# trunk/grimoire/immune-sentinel-memory.d/init.d.on/sentinel-memory
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: Prevents memory depletion
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

# Modifications to this file are LOST!
# Set variables SDIR and SIZE to alternate values only in
# /etc/init.d/conf.d/sentinel-memory

SDIR=/var/cache/swap		# directory which will hold swap files
SIZE=256M			# size of swap file

. /lib/lsb/init-functions

only start stop restart status
server	/usr/sbin/sentinel-memory
options	"$SDIR $SIZE"
