# trunk/grimoire/httpd.d/init.d/httpd
#!/bin/bash

### BEGIN INIT INFO
# Provides: httpd
# Required-Start: $network
# Should-Start: $remote_fs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: httpd is a feature packed web server
# Description: httpd is the most widely used web server.
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/httpd.pid
SERV=/usr/sbin/apachectl
NAME=httpd

status()  {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}

start()  {
  log_warning_msg  "$NAME about to start"
  [      -d  /var/log/httpd  ]  ||
  mkdir  -p  /var/log/httpd  ]
# if    start_daemon     -p $PIDF $SERV startssl
  if    start_daemon              $SERV startssl
  then  log_success_msg  "$NAME started"
  else  log_failure_msg  "$NAME not started"
        status
  fi
}

stop()  {
  log_warning_msg "$NAME about to stop"
  $SERV stop
  log_success_msg "$NAME stopped"
#  if    killproc        -p $PIDF -s 15  $SERV
#  then  log_success_msg "$NAME stopped"
#  else  log_failure_msg "$NAME not stopped"
#        status
#        return 0
#  fi
}

restart() {  $SERV restart;  }

try_restart() {
  status  &&  restart
  return 0
}

# reload()  {
#   PIDS="$( pidofproc -p $PIDF $SERV )"
#   if    [ -n            "$PIDS" ]
#   then  kill  -s SIGHUP  $PIDS  &&
#         log_success_msg "$NAME configuration reloaded"
#   fi
# }

help()  {
  echo "Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}

case  $1  in
       restart)      restart  ;;
   try-restart)  try_restart  ;;
        reload)  reload       ;;
  force-reload)  reload       ;;
         start)  start        ;;
        status)  status       ;;
          stop)  stop         ;;
             *)  help         ;;
esac
