# trunk/grimoire/sorcery-immune.d/init.d.on/immune
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Should-Start: $network
# Should-Stop:  $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: immune an immune system
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/immune.pid
SERV=/usr/sbin/immune
OPTS='critical neglect stall'

parse "$@"
