# trunk/grimoire/net-snmp.d/init.d/snmpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: $network run
# Required-Stop: $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: snmpd starts and stops the net-snmp SNMP daemon
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/snmpd
PIDF=/var/run/snmpd.pid
OPTS='-r -u nobody -p /var/run/snmpd.pid'

parse "$@"
