# trunk/grimoire/bsdftpd-ssl.d/xinetd.d/ftpd
service ftp
{
	flags		= REUSE NAMEINARGS
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/sbin/ftpd
	server_args	= ftpd  -l  -T 300  -z nosecure
	instances	= 128
	per_source	= 4
}
