# trunk/grimoire/pcmcia-cs.d/init.d/pcmcia-cs
#!/bin/bash

### BEGIN INIT INFO
# Provides: pcmcia-cs
# Required-Start: $syslog
# Required-Stop: $syslog
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: pcmcia-cs starts card services daemon
# Description: pcmcia-cs starts card services daemon
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg  "pcmcia-cs about to start"
  /etc/init.d/pcmcia start
  log_success_msg  "pcmcia-cs started"
}

stop()  {
  log_warning_msg  "pcmcia-cs about to stop"
  /etc/init.d/pcmcia stop
  log_success_msg  "pcmcia-cs stopped"
}

help()  {
  echo "Usage: $0 {start|stop}" 1>&2
}

case  $1  in
  start)  start  ;;
   stop)  stop   ;;
      *)  help   ;;
esac
