# trunk/grimoire/nsd.d/init.d.on/nsdc
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: $network run
# Required-Stop: $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: nsdc provides the name server daemon control script
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/nsd.pid
SERV=/usr/sbin/nsdc
CONF=/etc/nsd.conf
OPTS=$1

parse "$@"
