# trunk/grimoire/nsd.d/init.d.on/zonec
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Default-Start: 2 3 4 5
# Short-Description: zonec compiles a database for nsd
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/zonec
CONF=/etc/zonec/zonec.conf

parse "$@"
