# trunk/grimoire/immune-boot-kick.d/init.d/boot-kick-stop
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: S
# Short-Description: stops the boot-kick program before starting runlevels
### END INIT INFO

. /lib/lsb/init-functions

only start stop

start(){
 if   [ -x /etc/init.d/boot-kick ]
 then      /etc/init.d/boot-kick stop
 fi
 setleds -L < /dev/console > /dev/null
}
