# trunk/grimoire/deployment-docent.d/init.d/docent-deploy
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network empty run $local_fs
# Should-Start: dbus-daemon var_tmp
# Default-Start: 1 2 3 4 5
# Short-Description: a method for logging into the rootfs via ssh
### END INIT INFO

## /etc/init.d/deport-alien start
## causes creation of /etc/sorcery/local.config
# [ -f /etc/sorcery/local.config ] && exit 0

/bin/grep -q "FEATURE=()" /etc/sorcery/local.config && exit 0

. /lib/lsb/init-functions

only start stop
server /usr/sbin/docent-deploy
