# trunk/grimoire/libselinux.d/init.d.on/selinuxfs
#!/bin/bash

### BEGIN INIT INFO
# Provides: selinuxfs
# Required-Start: proc
# Should-Start:
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: selinuxfs mounts selinuxfs on /selinux
### END INIT INFO

. /lib/lsb/init-functions

NAME='/selinux'
OPTS='-o defaults,size=1m,nr_inodes=1k,mode=0755 -t selinuxfs selinuxfs'

parse "$@"
