# trunk/grimoire/immune-sentinel-space.d/init.d.on/sentinel-space
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: monitors and creates free space on file systems
### END INIT INFO

. /lib/lsb/init-functions

pid_file /var/run/sentinel-space.pid
server   /sbin/sentinel-space

parse "$@"
