# trunk/grimoire/immune-sentinel-space.d/init.d.on/sentinel-space
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: creates free space
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

. /lib/lsb/init-functions

server /usr/sbin/sentinel-space
