# trunk/grimoire/agrsm.d/init.d.on/agrserial
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs agrmodem
# Should-Start: udev
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Short-Description: agrserial loads the module
### END INIT INFO

. /lib/lsb/init-functions
