# trunk/grimoire/agrsm.d/init.d.on/agrserial
#!/bin/bash

### BEGIN INIT INFO
# Provides: agrserial
# Required-Start: $local_fs agrmodem
# Should-Start: udev
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Default-Stop:
# Short-Description: agrserial loads the module
# Description: agrserial loads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME=agrserial

parse "$@"
