# trunk/grimoire/LVM2.d/init.d/mapper
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dmsetup
# Should-Start: udev mdadm
# Default-Start: S
# Short-Description: loads map files from /etc/mapper into device-mapper
### END INIT INFO

. /lib/lsb/init-functions

cgroup none
only start
config /etc/mapper

start(){
run(){
 for MAP in /etc/mapper/*; do
  [ -f $MAP ] || continue
  log_warning_msg "mapper executing dmsetup create ${MAP##*/} $MAP"
  dmsetup create ${MAP##*/} $MAP
 done
}

 [ -d /etc/mapper ] &&
 if   log_warning_msg 'mapper starting'; run
 then log_success_msg 'mapper started'
 else log_failure_msg 'mapper not started'
 fi
}
