# trunk/grimoire/LVM2.d/init.d/vgchange
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dmsetup lock sys
# Should-Start: mdadm udev multipathd vgscan
# Required-Stop: finale halt mdadm--assemble--scan mdadm--auto-detect reboot reboot-kexec
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: vgchange activates logical volumes
### END INIT INFO

rootfs || exit 0

. /lib/lsb/init-functions

cgroup none
only start stop
name vgchange
options '--ignorelockingfailure -ay'
server /sbin/vgchange

start(){
 if   log_warning_msg "$NAME $OPTS"; $SERV $OPTS
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

stop(){
 OPTS='--ignorelockingfailure -an'
 [ -f /proc/mounts ] &&
 if   log_warning_msg "$NAME $OPTS"; $SERV $OPTS
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}
