# trunk/grimoire/sorcerer-sentient.d/init.d.on/sentient
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: gui_login
# Short-Description: Starts kdm-v4 after sorcerer-sentient installs
### END INIT INFO

[ "$1" == try-restart ] || exit 2

/etc/init.d/kdm-v4 status &>/dev/null ||
/etc/init.d/kdm-v4 start
