# trunk/grimoire/sorcerer-sentient.d/init.d.on/sentient
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: gui_login
# Short-Description: Reboots after sorcerer-sentient installs
### END INIT INFO

. /lib/lsb/init-functions
only try-restart

try_restart(){
 augur queue dispel/immune-boot-sentient
 dispel
 /etc/init.d/kdm status &>/dev/null ||
 shutdown -r +1
}
