# trunk/grimoire/sorcerer-sentient.d/init.d.on/sentient
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: gui_login
# Short-Description: Reboots after sorcerer-sentient installs
### END INIT INFO

. /lib/lsb/init-functions
only try-restart

try_restart(){
 dispel_restart(){
  while pgrep -u root "^cast$"
  do    sleep 30
  done
  if [ -f /etc/sorcery/log/version/immune-boot-sentient ]; then
   augur queue dispel/immune-boot-sentient
   dispel
   while pgrep -u root "^dispel$"
   do    sleep 15
   done
  fi
  if [ -f /root/restart-sentient ]; then
   rm  -f /root/restart-sentient ]
   shutdown -r now
  fi
 }

 dispel_restart &
}
