# trunk/grimoire/X.org.d/init.d/ICEperms
#!/bin/sh

### BEGIN INIT INFO
# Provides: bootmisc
# Required-Start: $local_fs
# Should-Start:
# Required-Stop: 
# Default-Start: S
# Default-Stop: 
# Short-Description: bootmisc performs miscellaneous actions on boot
# Description: bootmisc performs miscellaneous actions on boot
### END INIT INFO

. /lib/lsb/init-functions

log_success_msg "creating ICE-unix and setting perms"

mkdir  -p  /tmp/.ICE-unix
chown root /tmp/.ICE-unix
chmod 1777 /tmp/.ICE-unix

