# trunk/grimoire/openldap.d/init.d/slapd
#!/bin/bash

### BEGIN INIT INFO
# Provides: slap
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: $network
# Should-Stop:  $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: slapd is a stand-alone LDAP directory server
# Description: slapd is a stand-alone LDAP directory server
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/slapd.pid
SERV=/usr/libexec/slapd
NAME=slapd

parse "$@"
