# trunk/grimoire/init-scripts.d/init.d/bind
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dev-bind root-rw root-rw-retry root-alt
# Default-Start: S
# Short-Description: binds mounted virtual filesystems into /media/root/*/
#	except cgroup and dev
### END INIT INFO

. /lib/lsb/init-functions
only start
deny control

start(){

 vfs(){
  cut -d ' ' -f2 /proc/mounts |
  sed '\,^/$,d;\,^/media/root,d;\,^/dev$,d;\,/cgroup,d' |
  LC_ALL=C sort
 }

 b(){ [ -d "$2" ] || mkdir -pm 755 "$2"; mount --bind "$1" "$2"; }

 bind(){
   if   log_warning_msg "mount --bind $1 $2 executing"; b "$1" "$2"
   then log_success_msg "mount --bind $1 $2 executed"
   else log_failure_msg "mount --bind $1 $2 failed"; return 1
   fi
 }

 run(){
  local r=0
  for   f in $( vfs ); do
   for  d in /media/root/*
   do   bind "$f" "$d$f" || r=1
   done
  done
  return $r
 }

 has_root(){
  for d in /media/root/*/usr
  do  [ -d $d ] && return 0
  done
  exit 0
 }

 has_root
 if   log_warning_msg 'Virtual filesystems to /media/root/*/ binding'; run
 then log_success_msg 'Virtual filesystems to /media/root/*/ bound'
 else log_failure_msg 'Virtual filesystems to /media/root/*/ bound failed'
 fi
}
