# trunk/grimoire/init-scripts.d/init.d/bootlogd
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# No optimal time for stopping.
# Required-Start: pts ptmx
# Required-Stop: $local_fs kill
# Short-Description: runs console message capture logging daemon
### END INIT INFO

. /lib/lsb/init-functions

only start stop status
server /sbin/bootlogd

if   [ -d /home ] || grep -q boot=IR /proc/cmdline
then options -l /var/log/console
else options -l /media/root/var/log/console
fi
