# trunk/grimoire/init-scripts.d/init.d/btmp
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: creates login record /var/log/btmp
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control
name /var/log/btmp

start(){
 run(){
  if   ! [   -d ${NAME%/*} ]
  then mkdir -p ${NAME%/*}
  fi &&
  echo  -n  >    $NAME &&
  chown root:tty $NAME &&
  chmod 600      $NAME
 }

 if   log_warning_msg "$NAME create"; run
 then log_success_msg "$NAME create"
 else log_failure_msg "$NAME create"; return 1
 fi
}
