# trunk/grimoire/init-scripts.d/init.d/cgroup_root
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc sys
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts tmpfs file system on /sys/fs/cgroup
# Long-Description: Read /usr/src/linux/Documentation/cgroups/cgroups.txt
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /sys/fs/cgroup
options '-o defaults,size=64k,nr_inodes=16,mode=755 -t tmpfs cgroup_root'
