# trunk/grimoire/init-scripts.d/init.d/cgroup_cpuset
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc sys cgroup_root
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts cgroup file system on /sys/fs/cgroup/cpuset
# Long-Description: Read /usr/src/linux/Documentation/cgroups/cpusets.txt
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /sys/fs/cgroup/cpuset
options '-t cgroup -o cpuset,clone_children cpuset'
