# trunk/grimoire/init-scripts.d/init.d/cgroup_cpu
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: cgroup_root
# Default-Start: S
# Short-Description: mounts cgroup file system on /sys/fs/cgroup/cpu
# Long-Description: Read /usr/src/linux/Documentation/cgroups/cgroups.txt
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /sys/fs/cgroup/cpu
options '-t cgroup -o cpu,clone_children cpu'
