# trunk/grimoire/init-scripts.d/init.d/clean_tmp
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: removes /aux/run/{pipe,save,source,tmp}/* /tmp/*
### END INIT INFO

(( UID == 0 )) || exit 0

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  rm   -fr /aux/run/{pipe,save,source,tmp}/*
  if    ! [       -d      /tmp/ ]
  then /bin/rm    -f      /tmp/
       /bin/mkdir -m 1777 /tmp/
  else /bin/find          /tmp/ -mindepth 1 -maxdepth 1 -print0 |
       /bin/xargs -r -0   --max-lines=64 \
       /bin/rm -fr
  fi
  # Correct ownership and mode for /tmp set by /etc/init.d/var_tmp
 }

 if   log_warning_msg "/tmp/ cleaning"; run
 then log_success_msg "/tmp/ cleaned"
 else log_failure_msg "/tmp/ cleaning failed"
 fi
}
