# trunk/grimoire/init-scripts.d/init.d/clean_tmp
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: removes /fastboot and /forcefsck /tmp/* /var/cache/sorcery/*
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  renew(){
   if   [ -d $1 ]
   then find $1 -mindepth 1 -maxdepth 1 |
        xargs -r --max-lines=64 rm -fr
        chmod $2 $1
   else mkdir -pm $2 $1
   fi
  } 
  if   [ -f /fastboot  ]; then rm -f /fastboot;  touch /init; fi
  if   [ -f /forcefsck ]; then rm -f /forcefsck; touch /init; fi
  renew /tmp 1777
  if   [ -d /var/cache/sorcery ]
  then find /var/cache/sorcery -mindepth 1 -maxdepth 1 -type d | sed '\,/ccache,d'| xargs -r --max-lines=64 rm -fr
  fi
 }

 if   log_warning_msg "$NAME running"; run
 then log_success_msg "$NAME ran"
 else log_failure_msg "$NAME errored"
 fi
}
