# trunk/grimoire/init-scripts.d/init.d/config.gz
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: proc
# Default-Start: S
# Short-Description: chmod 600 /proc/config.gz
### END INIT INFO

. /lib/lsb/init-functions

only start
name /proc/config.gz

start(){
 [ -f $NAME ] || return
 if   log_warning_msg "$NAME changing mode to 600"; chmod 600 $NAME
 then log_success_msg "$NAME mode changed"
 else log_failure_msg "$NAME mode not changed"
 fi
}
