# trunk/grimoire/init-scripts.d/init.d/config
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Required-Start: sys
# Short-Description: mounts configfs on /config
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /config
options '-n -t configfs configfs'
