# trunk/grimoire/init-scripts.d/init.d/console
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 1 2 3 4 5
# Short-Description: creates /var/log/console
### END INIT INFO

. /lib/lsb/init-functions

only start
name /var/log/console

start(){ 
 run(){
  if   log_warning_msg "create"; echo -n > "/var/log/console"
  then log_success_msg "create"
  else log_failure_msg "create"; return 1
  fi
 }

 [ -f /var/log/console ] || run
}
