# trunk/grimoire/init-scripts.d/init.d/cpufreq_ondemand_ignore_nice_load
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: cpufreq_ondemand
# Default-Start: S
# Short-Description: ondemand will not elevate CPU freq for nice > 0 processes
### END INIT INFO

. /lib/lsb/init-functions
only start stop status

inl=/sys/devices/system/cpu/cpufreq/ondemand/ignore_nice_load

start(){  if [ -f $inl ]; then echo 1 > $inl; fi; }
stop(){   if [ -f $inl ]; then echo 0 > $inl; fi; }
status(){ if [ -f $inl ]; then cat      $inl; fi; }
