# trunk/grimoire/init-scripts.d/init.d/depmod
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: runs depmod during boot
### END INIT INFO

. /lib/lsb/init-functions

fast
cgroup none
only start
server /sbin/depmod
options -A
