# trunk/grimoire/init-scripts.d/init.d/empty
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts tmpfs on /var/empty
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /var/empty
options '-o defaults,size=64k,nr_inodes=16,mode=700 -t tmpfs tmpfs'
