# trunk/grimoire/init-scripts.d/init.d/etc-rootfs
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Short-Description: copies /aux/set/rootfs/etc/ to / on rootfs.
### END INIT INFO

. /lib/lsb/init-functions

if ! rootfs; then trap - EXIT; exit 0; fi

only start
deny control
config 

start(){
 copy(){
  /bin/cp -a /etc/hostname        /etc/origin
  /bin/cp -a /aux/set/rootfs/etc/ /
  return 0
 }

 if   log_warning_msg "/aux/set/rootfs/etc/ to / copy"; copy
 then log_success_msg "/aux/set/rootfs/etc/ to / copy"
 else log_failure_msg "/aux/set/rootfs/etc/ to / copy"; return 1
 fi
}
