# trunk/grimoire/init-scripts.d/init.d/final_check
#!/bin/bash
# Copyright 2011,2012,2013 by Kyle Sallee
# All Rights Reserved
# For use on Modern Magic only

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 0 6
# Short-Description: if critical error then abort reboot/shutdown and begin runlevel 3
### END INIT INFO

# To disable rebound
# then edit /etc/init.d/conf.d/halt
# and add function declarations:
# rebound(){ :; }

. /lib/lsb/init-functions

# do not execute on rootfs
# only execute on root file systems.
if rootfs; then trap - EXIT; exit 0; fi

only start
deny control

start(){

 rebound(){
  [ -f /proc/mounts ] || return
  /bin/find /var/log/init/ -mindepth 1 -maxdepth 1 |
  while read; do /bin/mv $REPLY{,.old}; done
  /etc/init.d/rcS
  /sbin/telinit $REPLY
  exit
 }

 critical(){
  alert(){
   short(){ echo -en "\a"; }
    long(){ for ((CX=0;CX<35;CX++)); do echo -en "\a"; sleep 0.0016; done; }
   pause(){ sleep 0.5; }

   S(){ echo -n "S"; short; pause; short; pause; short; pause; pause; }
   O(){ echo -n "O"; long;  pause; long;  pause; long;  pause; pause; }

   S; O; S; echo; sleep 1; S; O; S;
  }

  warn(){
   echo "/etc/init.d/$1 reported an error"
   echo "reboot/shutdown request denied"
   echo "Rebounded to runlevel 3"
   echo "This message originates from /etc/issue"
   echo "Please edit the file and remove the error message"
   echo "after the error condition is fixed"
  }

  deny(){ warn $1 | tee -a /etc/issue; alert; REPLY=3; rebound; }

    if [ -f /var/log/init/failed/extlinux/stop  ]; then deny extlinux
  elif [ -f /var/log/init/failed/initramfs/stop ]; then deny initramfs
  fi
 }

 if   log_warning_msg "Final check"; critical
 then log_success_msg "Final check"
 else log_failure_msg "Final check"; return 1
 fi
}
