# trunk/grimoire/init-scripts.d/init.d/finale
#!/bin/bash
# Copyright 2011, 2012, 2013 by Kyle Sallee
# All Rights Reserved
# For use on Sorcerer only

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 0 6
# Required-Start: final_check
# Short-Description: final pause, login, and rebound before the halt or reboot
### END INIT INFO

. /lib/lsb/init-functions

if ! rootfs; then trap - EXIT; exit 0; fi

only start
deny control

start(){

rebound(){
 [ -f /proc/mounts ] || return
 /bin/find /var/log/init/ -mindepth 1 -maxdepth 1 |
 while read; do /bin/mv $REPLY{,.old}; done
 /etc/init.d/rcS
 /sbin/telinit $REPLY
 exit
}

count_down(){
 local CX=10
 echo -n "(p)ause; (s)ulogin; (#)runlevel"
 while (( CX-- > 0 )); do
  if read -p " $CX" -n 1 -t 1; then
   echo
   case $REPLY in
    p) read -n 1 -t 300 ;;
    s) sulogin   -t 300 ;;
    [0-7]) rebound      ;;
    *) CX=0             ;;
   esac
  fi
 done
}

 if   log_warning_msg "Finale begins"; count_down
 then log_success_msg "Finale complete"
 else log_failure_msg "Finale error"
 fi
}
