# trunk/grimoire/init-scripts.d/init.d/fstab
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: alt_boot find_ir udevd
# Default-Start: S
# Short-Description: avoids rdinit=/bin/bash when / is not in fstab
# Description: fstab checks for / in /etc/fstab and rewrites it to /media/root/real
### END INIT INFO

if   [ -d /home ] || grep -q boot=IR /proc/cmdline
then exit 0
fi

. /lib/lsb/init-functions

only start
deny control

start(){
 root_menu(){
  block(){ find /dev -type b -printf '%p	%Y\n' | sed '/loop/d'; }

  dialog --backtitle fstab --stdout --no-cancel --title "Root Device Menu" \
         --menu "Expected / in /etc/fstab" 0 0 0 /sbin/sulogin '' $( block )
 }

 interact(){
  # usb disk device may not be immediately available?
  sleep 5
  ROOT=$( root_menu )
  case $ROOT in
   /sbin/sulogin) sulogin -t 300 1>/dev/console 2>&1 0</dev/console ;;
   /dev*) echo "$ROOT /media/root/real auto defaults 0 0" >> /etc/fstab.rootfs
  esac
 }

 check(){ [ -n "$( sed "s/#.*//;s:\t: :g;\: $1:p;d" /etc/fstab.rootfs )" ]; }

 copy(){
  sed 's:#.*::' /etc/fstab |
  sed 's:\t: :g' |
  sed 's: / : /media/root/real :p;d' > /etc/fstab.rootfs
 }

 [ -f /etc/fstab.rootfs ] || copy

 if   log_warning_msg '/etc/fstab contains /media/root checking'; check /media/root/
 then log_success_msg '/etc/fstab contains /media/root verified'
 else log_failure_msg 'root file system search failed'; interact
 fi
}
