# trunk/grimoire/init-scripts.d/init.d/halt
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 0
# Required-Start: finale
# Short-Description: powers down the box
### END INIT INFO

. /lib/lsb/init-functions

if ! rootfs; then trap - EXIT; exit 0; fi

only start
deny control

start(){ 
 log_success_msg Halting
 /sbin/halt -d -f -i -p
}
