# trunk/grimoire/init-scripts.d/init.d/hostname
#!/bin/bash

# Copyright Kyle Sallee 2011, 2012, 2013
# All rights reserved.
# For use with the distribution Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S
# Short-Description: sets the hostname for the rootfs
# run /etc/init.d/rootname configure to set a hostname for a root file system.
### END INIT INFO

. /lib/lsb/init-functions

only start configure

configure(){
 hostname_dialog(){
  BACKTITLE=$"hostname configuration menu"
      TITLE=$"hostname input"
       HELP=$"
Rootfs and root file systems each receive a FQDN,\n
fully qualified domain names.\n
When a computer serves multiple roles\n
then the hostname for the rootfs could be the domain name.\n
Then the hostname for each root file system\n
could be the root file system's name dot hostname.\n
If the hostname for rootfs is party.tea\n
then /+/base could could be named base.party.tea,\n
and  /+/www  could could be named  www.party.tea.\n
Hostnames must be unique.\n
Hostnames for root file systems\n
can entirely differ from the hostname for the rootfs.\n
\n
The hostname for the rootfs        is written to /etc/rootfs/hostname\n
The hostname for root file systems is written to /etc/hostname\n
The init-script for configuring the hostname\n
for a root file system is called /etc/rootname"

  /bin/dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout --no-cancel      \
   --inputbox "$HELP" 0 0
 }

 confirm_it(){ /bin/dialog --nocancel --yesno "Missing dot.  Use $NAME as is?" 0 0; }

 /bin/mkdir -pm 700  /etc/rootfs/
 local NAME=$( hostname_dialog )
 if  echo "$NAME" | /bin/grep -q "\." || confirm_it; then
     echo "$NAME" >        /etc/rootfs/hostname
    if   [ -d       /+/base ]; then
     mkdir -pm 700  /+/base/etc/rootfs/
     echo "$NAME" > /+/base/etc/rootfs/hostname
    fi
 fi
}

start(){
 rootfs || exit 0
 local REPLY
 log_warning_msg "hostname setting"
 if   [ -f /etc/hostname ]; then read REPLY < /etc/hostname; fi
 if   [ -z "$REPLY" ]
 then REPLY="unknown$RANDOM$RANDOM"
      log_failure_msg "/etc/init.d/hostname configure"
      log_failure_msg "Executing the above command is recommended"
 fi

 echo "$REPLY" > /proc/sys/kernel/hostname
 log_success_msg "hostname is now $REPLY"
}
