# trunk/grimoire/init-scripts.d/init.d/hwclock
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: hwclock localtime
# Required-Start: udev
# Should-Start: rtc-cmos
# Required-Stop:  $local_fs
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: loads and saves system time; hardware clock uses UTC
### END INIT INFO

# Either something resets system to hardware clock UTC during sysinit
# or restoring system time from hwclock during sysinit
# does not work as expected.
# Therefore starting at S 1 2 3 4 5

. /lib/lsb/init-functions

NAME='System clock from hardware clock time restoring'
SERV=/sbin/hwclock

start(){
 if   log_warning_msg 'System clock from hardware clock'; /sbin/hwclock -s
 then log_success_msg 'System clock set'
 else log_failure_msg 'System clock not set'
 fi
}

stop(){
 if   log_warning_msg 'Hardware clock from system clock'; /sbin/hwclock -w
 then log_success_msg 'Hardware clock set'
 else log_failure_msg 'Hardware clock not set'
 fi
}

restart(){ :; }

parse "$@"
