# trunk/grimoire/init-scripts.d/init.d/hwclock-localtime
#!/bin/bash

### BEGIN INIT INFO
# Provides: hwclock localtime
# Required-Start: udev
# Should-Start: rtc-cmos
# Required-Stop:  $local_fs
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: loads and saves system time; hardware clock uses local time
# When started in sysinit something, perhaps fsck, restores time from
# the hardware clock in UTC format rendering hwclock-locatime ineffective
# That is why hwclock-localtime is started both in sysinit and runlevels.
### END INIT INFO

. /lib/lsb/init-functions

NAME='System clock from hardware clock time restoring'
SERV=/sbin/hwclock
OPTS='--hctosys --localtime'

stop(){
 if   log_warning_msg 'Hardware clock from system time setting'; /sbin/hwclock --systohc --noadjfile --localtime
 then log_success_msg 'Hardware clock set'
 else log_failure_msg 'Hardware clock not set'
 fi
}

restart(){ :; }

parse "$@"
