# trunk/grimoire/init-scripts.d/init.d/hwclock
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 400
# Provides: hwclock localtime
# Required-Start: udev
# Should-Start: rtc-cmos
# Required-Stop:  $local_fs
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: loads and saves system time; hardware clock uses UTC
### END INIT INFO

# Either something resets system to hardware clock UTC during sysinit
# or restoring system time from hwclock during sysinit
# does not work as expected.
# Therefore starting at S 1 2 3 4 5

. /lib/lsb/init-functions

only start stop restart configure display
server /sbin/hwclock

display(){ echo "Mode $MODE"; }

configure(){
 mode_menu(){
  BACKTITLE=$"hwclock configuration menu"
      TITLE=$"Mode Selection"
       HELP=$"POSIX/Unix/Linux boxes store time internally using UTC.
Do not select localtime unless certain it is required.
Please select either UTC or localtime."

  dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout --no-cancel      \
   --menu "$HELP" 0 0 0      \
   UTC Normal localtime Strange
 }

 case $( mode_menu ) in
  UTC)       echo "MODE=--utc" ;;
  localtime) echo "MODE=--localtime" ;;
 esac > /etc/init.d/conf.d/hwclock
}

start(){
 if   log_warning_msg 'System clock set'; /sbin/hwclock -s ${MODE:---utc}
 then log_success_msg 'System clock set'
 else log_failure_msg 'System clock set'; return 1
 fi
}

stop(){
 if   log_warning_msg 'Hardware clock set'; /sbin/hwclock -w ${MODE:---utc}
 then log_success_msg 'Hardware clock set'
 else log_failure_msg 'Hardware clock set'; return 1
 fi
}

restart(){ :; }
