# trunk/grimoire/init-scripts.d/init.d/ip_forward
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: proc
# Default-Start: S
# Short-Description: enables IP forwarding, necessary for veth
### END INIT INFO

. /lib/lsb/init-functions

only start stop status
config /proc/sys/net/ipv4/ip_forward

start(){
 if   log_warning_msg "ip_forward enabling"; echo 1 > /proc/sys/net/ipv4/ip_forward
 then log_success_msg "ip_forward enabled"
 else log_failure_msg "ip_forward enable failed"; return 1
 fi
}

stop(){
 if   log_warning_msg "ip_forward stopping"; echo 0 > /proc/sys/net/ipv4/ip_forward
 then log_success_msg "ip_forward stoped"
 else log_failure_msg "ip_forward stop failed"; return 1
 fi
}
