# trunk/grimoire/init-scripts.d/init.d/lang
#!/bin/bash

# Copyright Kyle Sallee 2011, 2012, 2013
# All rights reserved.
# For use with the distribution Sorcerer Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: S
# Short-Description: checks and sets LANG in /etc/bash.d/noninteractive/LANG
### END INIT INFO

# omit log

e=
for  d in /media/root/*; do
 if   [ -f   "$d"/usr/bin/locale ]
 then        "$d" $0 "$@"; e=1
 fi
done
if [ -n "$e" ]; then exit 0; fi

. /lib/lsb/init-functions

only start configure status
name /etc/bash.d/noninteractive/LANG

if   [ -f   /usr/lib64/locale/locale-archive ] 
then config /usr/lib64/locale/locale-archive 
else config   /usr/lib/locale/locale-archive 
fi

configure_run(){

 lang_list(){ locale -a | sed "s:$:	locale:"; }

 lang_menu(){
  BACKTITLE=$"Locale Configuration Menu"
      TITLE=$"Locale Selection"
       HELP=$"Please select the computer's default locale within 60 seconds."

  dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout --no-cancel      \
   --menu "$HELP" 0 0 0      \
   $( lang_list )
 }

 LANG="$( lang_menu )"
 if [ "$LANG" ]; then
  rm -f  $NAME
  mkdir -pm 755 ${NAME%/*}
  echo "export LANG=$LANG" > $NAME
  chmod 755 $NAME
 fi
}

status(){
 if   log_warning_msg "$NAME checking"; [ -f $NAME ]
 then log_success_msg "$NAME checked"
 else log_warning_msg "$NAME missing"; return 1
 fi
}

choices(){ (( $( locale -a 2>/dev/null | wc -l ) > 2 )); }
extreme(){ (( $( locale -a 2>/dev/null | wc -l ) > 9 )); }

trim_locale_archive(){
 if   extreme &&
      dialog --title "Want to trim locale-archive of unsused locales?" \
             --defaultno --yesno \
             "Run /etc/init.d/locale-archive configure?" 0 0
 then /etc/init.d/locale-archive configure
 fi
}

make_locale_archive(){
 if ! choices; then
  /etc/init.d/locale-archive configure
  echo "Please be patient while locale-archive is generated."
  if   [      -f /usr/libexec/bash/sleep ]
  then enable -f /usr/libexec/bash/sleep sleep
  fi
  while (( $( pgrep -u root locale-archive | wc -l ) > 1 )); do sleep 1; done
 fi
}

start(){ if choices && ! status; then configure_run; fi; }
configure(){
 rm -f $NAME
 make_locale_archive
 configure_run
 trim_locale_archive
}
