# trunk/grimoire/init-scripts.d/init.d/locale-archive
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: 2 3 4 5
# Short-Description: generates locales databases
### END INIT INFO

. /lib/lsb/init-functions

if rootfs; then trap - EXIT; exit 0; fi

only start configure status reload display
cgroup idle
config /usr/share/i18n/SUPPORTED
server /sbin/locale-archive

display(){ echo "Locales $LOCALES"; }

status(){
 local LA
 if   [[ $HOSTTYPE == x86_64 ]]
 then LA=/usr/lib64/locale/locale-archive
 else LA=/usr/lib/locale/locale-archive
 fi

 if   log_warning_msg "$LA checking"; ! [ -f $LA ]
 then log_warning_msg "$LA missing"; return 1
 elif [ /usr/bin/localedef -nt $LA ]
 then log_warning_msg "$LA obsolete"; return 1
 else log_success_msg "$LA checked"
 fi
}

choices(){
 local CX=$( locale -a 2>/dev/null | wc -l )
 (( CX > 2 ))
}

configure(){
 /bin/tty -s || return

 locales_checklist(){
  sed -n 's: .*::;/\//p' /usr/share/i18n/SUPPORTED |
  while read; do
     NAME="${REPLY%/*}"
   SOURCE="${NAME%.*}"

   echo "$REPLY"
   sed  -nr 's:title[^"]*"([^"]*)".*:\1:p' /usr/share/i18n/locales/$SOURCE
   echo "off"
  done
 }

 locales_menulist(){ 
  if   [   -n "$LOCALES" ]
  then echo   "$LOCALES" |
       sed -n 's: .*::;/\//p'
  else sed -n 's: .*::;/\//p' /usr/share/i18n/SUPPORTED
  fi |
  while read; do
     NAME="${REPLY%/*}"
    SOURCE="${NAME%.*}"

   echo "$NAME"
   sed -nr 's:^title[^"]*"([^"]*)".*:\1:p' /usr/share/i18n/locales/$SOURCE
  done
 }

 select_locales(){
  BACKTITLE=$"Glibc Locale Configuration Menu"
      TITLE=$"Locale Selection"
       HELP=$"Please select available language locales."

  dialog --backtitle  "$BACKTITLE" \
         --title      "$TITLE"     \
         --stdout                  \
         --separate-output         \
         --timeout  60             \
         --no-cancel               \
         --checklist  "$HELP"      \
         0 0 0                     \
         $( locales_checklist )
 }

 local  IFS
 export IFS=$'	\n'
 if   dialog --title "locale-archive configuration" \
             --defaultno --yesno \
             "Select individual locales?" 0 0
 then       LOCALES="$( select_locales )"
 else unset LOCALES
 fi
 echo "LOCALES=\"$LOCALES\"" > /etc/init.d/conf.d/locale-archive
 reload
}

reload(){
 rm -f /usr/lib{,64}/locale/locale-archive
 /etc/init.d/locale-archive start
}

if [ "$1" == start ] && status &>/dev/null; then start(){ :; }; fi
