# trunk/grimoire/init-scripts.d/init.d/Menu
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 500
# Short-Description: mass init-script auditing/control menu
### END INIT INFO

. /lib/lsb/init-functions

only configure
deny control

configure(){

init_commands(){
 local  LC_ALL=C
 export LC_ALL

 grep -ilr "BEGIN INIT INFO" /etc/init.d |
 xargs --max-lines=256 cat |
 sed 's/^only//p;d' |
 tr '\t' '\n' |
 tr ' '  '\n' |
 sed '/^$/d'  |
 sort -u      |
 sed 's:$:\t Command:'
}

command_menu(){
 TITLE=$"Init-Script Selection"
 HELP=$"Please select an init-script command"
  dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout                  \
   --cancel-label "Cancel"   \
   --menu "$HELP" 0 0 0      \
   $( init_commands )
}


 init_scripts(){ grep -ilr "BEGIN INIT INFO" /etc/init.d | grep -vx "/etc/init.d/Menu" | sort; }

 add_mode(){
  while read; do find $REPLY -maxdepth 0 -printf "%p\t%m\n"; done |
  sed "s:\t510$:\tAuto Updates:
       s:\t500$:\tAuto        :
       s:\t400$:\tManual      :
       s:\t...$:\tCustom      :"
 }

 add_when(){
  while read NAME MORE; do
   START=$( sed "s/# Default-Start://p;d" $NAME | tr -s ' ' | sed "s/ / S/g" )
    STOP=$( sed "s/# Default-Stop://p;d"  $NAME | tr -s ' ' | sed "s/ / K/g" )
   echo "$NAME	$MORE $START$STOP"
  done
 }

 add_configure(){
  while read NAME MORE; do
   if   grep -q "only .*configure" $NAME
   then echo "$NAME	C $MORE"
   else echo "$NAME	  $MORE"
   fi
  done
 }

 add_item_help(){
  while read NAME REST; do
   echo "$NAME"
   echo "$REST"
   echo "off"
   DESC="$( sed 's/.*Short-Description: //p;d' /etc/init.d/$NAME )"
   echo "${DESC:-no short description}"
  done
 }

 init_list(){
  export IFS=" 	
"
  init_scripts |
  add_mode |
  add_configure |
  add_when |
  sed 's:.*/::' |
  add_item_help
 }

choices(){
 /etc/init.d/$1 help 2>&1 |
 sed "s:.*{::; s:|:\n:g; s:}::" |
 sed "s:$:\tChoice:"
}

confirm(){ dialog --cr-wrap --nocancel --yesno "$1" 0 0; }

init_script_checklist(){
 TITLE=$"Init-Script Selection"
 HELP=$"Please select init-scripts"

  local IFS="
"
  export IFS

  dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout --item-help      \
   --separate-output         \
   --ok-label "Done"         \
   --cancel-label "Abort"    \
   --checklist "$HELP" 0 0 0 \
   $( init_list )
 }

 BACKTITLE=$"Init-Script Menu"
 if CMD=$( command_menu ) && clear &&
     ISC=$( init_script_checklist ) &&
    ASK="$( echo "$ISC" | sed "s:^:/etc/init.d/:; s:$: $CMD:" )" &&
    confirm "Execute:\n$ASK"
 then
  clear
  for IS in $ISC
  do  bash /etc/init.d/$IS $CMD
  done
 fi
}
