# trunk/grimoire/init-scripts.d/init.d/modules
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: root-rw
# Default-Start: S
# Short-Description: provides modules to root file systems that lack.
### END INIT INFO

. /lib/lsb/init-functions

if   ! rootfs
then echo "$0 can only be executed on the rootfs"; trap - EXIT; exit 1
fi

only start
deny control

mod_copy(){
 local d="$1"
 if ! [  -d   "$d/lib/modules/$V/" ]
 then cp -a {,"$d"}/lib/modules/"$V"
 fi
}

mod_check(){ [ -d "$1/lib/modules/$V/" ]; }

start(){
 local V="$( cut -d ' ' -f3 /proc/version )"
 local r
 for d in /media/root/*; do
  if   ! [ -f "$d/init" ]; then continue; fi
  if   log_warning_msg "$d/lib/firmware/   checking"; ! [ -d "$d/lib/firmware" ]
  then log_success_msg "$d/lib/firmware/   missing"
       log_success_msg "$d/lib/modules/$V/ otiose"
  elif log_warning_msg "$d/lib/modules/    checking"; ! [ -d "$d/lib/modules" ]
  then log_success_msg "$d/lib/modules/    missing"
       log_success_msg "$d/lib/modules/$V/ otiose"
  elif log_warning_msg "$d/lib/modules/$V/ checking"; mod_check "$d"
  then log_success_msg "$d/lib/modules/$V/ exists"
  elif log_warning_msg "$d/lib/modules/$V/ copying";  mod_copy  "$d"
  then log_success_msg "$d/lib/modules/$V/ installed"
  else log_failure_msg "$d/lib/modules/$V/ copy failed"; r=1
  fi
 done
 return $r
}
