# trunk/grimoire/init-scripts.d/init.d/mount
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: $local_fs
# Required-Start: fsck mtab
# Should-Start: root-rw
# Required-Stop: cryptmount cryptsetup finale halt mdadm reboot reboot-kexec vgchange
# Should-Stop: root-ro
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts and un-mounts file systems
### END INIT INFO

. /lib/lsb/init-functions

only start stop
deny control

FSTAB=/etc/fstab

mkdirs(){
 while read s d t; do
  if   [[    $s =~ \#     ]]; then continue
  elif [[    $t =~ noauto ]]; then continue
  elif [[ -z $d           ]]; then continue; fi
  echo "$d"
 done < $FSTAB |
 /bin/xargs -r -d '\n' -L64 -- /bin/mkdir -p 0755 --
}

run(){ if [ -f $FSTAB ]; then mkdirs && /bin/mount -a; fi; }

start(){
 if   log_warning_msg "/bin/mount -a execute"; run
 then log_success_msg "/bin/mount -a execute"
 else log_success_msg "/bin/mount -a execute"
 fi
}

stop(){
 run(){
  sync
  [ -f /proc/mounts ] || mount -nrt proc proc /proc

  if   rootfs; then
   /bin/cut -d ' ' -f2 /etc/mtab /proc/mounts |
   /bin/sed '\,^/$,d' |
   /bin/sort -ru | /bin/xargs -r --max-lines=256 /bin/umount -l
  else
   /bin/cut -d ' ' -f2 /etc/mtab |
   /bin/sort -ru | /bin/xargs -r --max-lines=256 /bin/umount -l
  fi
 }

 log_warning_msg 'File systems unmounting'; run
 log_success_msg 'File systems unmounted'
}
