# trunk/grimoire/init-scripts.d/init.d/mount
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: $local_fs
# Required-Start: mtab root-rw dev-bind
# Required-Stop: cryptmount cryptsetup finale halt mdadm reboot reboot-kexec root-ro vgchange
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mounts and un-mounts file systems
### END INIT INFO

. /lib/lsb/init-functions

only start stop
deny control

start(){
 local r=0
 if    [ -d /home ]; then
  if   log_warning_msg "mount -a executing"; mount -a
  then log_success_msg "mount -a executed"
  else log_success_msg "mount -a execution failed"; r=1
  fi
 else
  for d in /media/root/*; do
   if [ "$d" == /media/root/rootfs ] || ! [[ -f $d/etc/fstab ]]
   then continue
   elif log_warning_msg "chroot $d mount -a executing"; chroot "$d" mount -a
   then log_success_msg "chroot $d mount -a executed"
   else log_success_msg "chroot $d mount -a execution failed"; r=1
   fi
  done
 fi
 return $r
}

stop(){
 run(){
  sync
  [ -f /proc/mounts ] || mount -nrt proc proc /proc

  cut  -d ' ' -f2 /proc/mounts  |
  sed  '\,^/media/root,p;d'     |
  sort -r | xargs -r --max-lines=256 umount -lr

  sed  '\,^rootfs / rootfs,d
        \,^proc /proc proc,d
        \,^sysfs /sys sysfs,d
        \,^tmpfs /run tmpfs,d
        \,^devtmpfs /dev devtmpfs,d
        \,^devpts /dev/pts devpts,d' /proc/mounts | cut -d ' ' -f2 |
  sort -r | xargs -r --max-lines=256 umount -lr
 }

 log_warning_msg 'File systems unmounting'; run
 log_success_msg 'File systems unmounted'
}
