# trunk/grimoire/init-scripts.d/init.d/mount
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: $local_fs
# Required-Start: mtab root-rw dev-bind
# Required-Stop: cryptmount cryptsetup finale halt mdadm reboot reboot-kexec root-ro vgchange
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mounts and un-mounts file systems
### END INIT INFO

. /lib/lsb/init-functions

only start stop
deny control

start(){
 run(){ chroot "$1" mount -a; }

 local r=0
 if    [ -d /home ]; then
  if   log_warning_msg "mount -a executing"; mount -a
  then log_success_msg "mount -a executed"
  else log_success_msg "mount -a execution failed"; r=1
  fi
 else
  for d in /media/root/*; do
   [[ -f $d/etc/fstab ]] || continue
   if   log_warning_msg "chroot $d mount -a executing"; run "$d"
   then log_success_msg "chroot $d mount -a executed"
   else log_success_msg "chroot $d mount -a execution failed"; r=1
   fi
  done
 fi
 return $r
}

stop(){
 run(){
  sync
  [ -f /proc/mounts ] || mount -nrt proc proc /proc

  cut -d ' ' -f2 /proc/mounts |
  LC_ALL=C sort -r |
  if   [ -d /media/root ]
  then grep /media/root
  else grep -vx / ; mount -o remount,ro /
  fi |
  xargs -r --max-lines=1024 umount -lr
 }

 log_warning_msg 'File systems unmounting'; run
 log_success_msg 'File systems unmounted'
}
