# trunk/grimoire/init-scripts.d/init.d/mqueue
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dev proc
# Should-Start: dev
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts mqueue on /dev/mqueue
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /dev/mqueue
options '-t mqueue mqueue'
