# trunk/grimoire/init-scripts.d/init.d/mtab
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: bind
# Default-Start: S 1 2 3 4 5
# Short-Description: sanatizes /etc/mtab
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  if   [ -d /etc/mtab ]; then /bin/rmdir /etc/mtab
  elif [ -e /etc/mtab ]; then /bin/rm -f /etc/mtab
  fi
  echo -n > /etc/mtab
 }

 if   log_warning_msg 'echo -n > /etc/mtab'; run
 then log_success_msg 'echo -n > /etc/mtab'
 else log_failure_msg 'echo -n > /etc/mtab'; return 1
 fi
}
