# trunk/grimoire/init-scripts.d/init.d/mtab
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: bind
# Default-Start: S
# Short-Description: creates /etc/mtab from /proc/mounts
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  rm -f                                   /etc/mtab
  ln -s                 /proc/mounts      /etc/mtab
  for d in   /media/root/*; do
   [ "$d" == /media/root/rootfs ] && return
   if   [ -f         "$d/proc/mounts" ]
   then chroot "$d" cat /proc/mounts > "$d/etc/mtab"
   fi
  done
 }

 log_warning_msg '/etc/mtab creating'; run
 log_success_msg '/etc/mtab created'
}
