# trunk/grimoire/init-scripts.d/init.d/mtab
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: root-rw
# Should-Start: root-rw-retry root-alt
# Default-Start: S
# Short-Description: creates /etc/mtab from /proc/mounts
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  MR=/media/root
  if   [  -f  $MR/sbin/init ]
  then rm -fr $MR/etc/mtab /etc/mtab; cat /proc/mounts > $MR/etc/mtab; ln -s /proc/mounts /etc/mtab
  else rm -fr              /etc/mtab; cat /proc/mounts >    /etc/mtab
  fi
 }

 log_warning_msg '/etc/mtab creating'; run
 log_success_msg '/etc/mtab created'
}
