# trunk/grimoire/init-scripts.d/init.d/mtab
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: bind
# Default-Start: S 1 2 3 4 5
# Short-Description: sanatizes /etc/mtab
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control

start(){
 run(){
  if   [  -h /etc/mtab ]; then rm -f /etc/mtab
  elif [  -d /etc/mtab ]; then rmdir /etc/mtab
  fi
  echo -n > /etc/mtab
 }

 log_warning_msg '/etc/mtab create'; run
 log_success_msg '/etc/mtab create'
}
