# trunk/grimoire/init-scripts.d/init.d/nr_hugepages
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

# Not required to reserve huge pages
# transparent huge page support is superior
# and does not require memory reservation.

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: proc
# Should-Start: dev hugetlbfs sys
# Default-Start: 1 2 3 4 5
# Short-Description: sets number of huge pages 
### END INIT INFO

NR=$( sed -r "s/MemTotal: *([^ ]*) kB/\1/p;d"     /proc/meminfo )
HS=$( sed -r "s/Hugepagesize: *([^ ]*) kB/\1/p;d" /proc/meminfo )
(( NR /= 10 ))
(( NR /= HS ))

. /lib/lsb/init-functions

fast
only start configure display

display(){ echo "NR=$NR"; }

configure(){

nr_dialog(){
 BACKTITLE=$"nr_hugepages configuration menu"
     TITLE=$"Instructions"
      HELP=$"Please provide the number of unswappable two megabyte huge pages to reserve"

 dialog \
  --backtitle "$BACKTITLE" \
  --title     "$TITLE"     \
  --timeout 60             \
  --stdout --no-cancel     \
  --default-item 3         \
  --inputbox "$HELP" 0 0 ${NR:-128}
}

 NR=$( nr_dialog )
 if   [ -n "$NR" ]
 then echo "NR=$NR" > /etc/init.d/conf.d/nr_hugepages
 fi
}

start(){
 run(){
  if   [ -f         /proc/sys/vm/nr_hugepages ]
  then echo "$NR" > /proc/sys/vm/nr_hugepages
  fi
 }

 if   log_warning_msg "$NAME setting"; run
 then log_success_msg "$NAME set"
 else log_failure_msg "$NAME not set"; return 1
 fi
}
