# trunk/grimoire/init-scripts.d/init.d/power-loss-recovery
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: sorcery
# Default-Start: 1 2 3 4 5
# Short-Description: restores box to working condition following power loss.
# Copyright 2012 by Kyle Sallee All Rights Reserved
# For use with the distribution Sorcerer only.
### END INIT INFO

. /lib/lsb/init-functions

only start

recover(){
 local  l r=0
 local  RATH="$PATH"
 local  PATH="$o/bin:$o/usr/bin:$RATH"
 export PATH

 find $i/ -type f -printf '%f\n' |
 while read; do
  l="$o/etc/sorcery/log/install/$REPLY"
  if ! [ -f "$l" ]; then r=1; continue; fi
  case $REPLY in
   glibc)    r=1; continue ;;
   ncurses)  r=1; continue ;;
   readline) r=1; continue ;;
  esac
  tar -C "$o" --no-recursion -cPT - < "$l" |
  tar -C /                   -pxf -
  rm  -f "$i/$REPLY"
 done
 /sbin/ldconfig &>/dev/null
 /usr/sbin/cast
 return $r
}

check(){
 if   [ -d "$i" ]
 then [ -z "$( find "$i/" -type f )" ]
 fi
}

recovery(){
 if   log_warning_msg 'files recovering'; recover
 then log_success_msg 'files recovered'
 else log_failure_msg 'file recovery failed'
 fi
}

start(){
 local i=/var/log/sorcery/queue/installing
 local o=/aux/old

 if   log_warning_msg 'power lost during cast checking'; check
 then log_success_msg 'power lost during cast checked';
 else recovery
 fi
}
