# trunk/grimoire/init-scripts.d/init.d/proc
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: mtab
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts proc on /proc
### END INIT INFO

. /lib/lsb/init-functions

only start restart status
deny control
name /proc
options '-n -t proc proc'
