# trunk/grimoire/init-scripts.d/init.d/pts
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Should-Start: dev
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts devpts on /dev/pts
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /dev/pts

# mi='-q CONFIG_DEVPTS_MULTIPLE_INSTANCES=y'
# ni='-t devpts -o newinstance -o ptmxmode=0666 devpts'
# c=/proc/config

# if   [ -f $c    ] &&                       /bin/grep $mi $c; then options "$ni"
# elif [ -f $c.gz ] && /bin/gzip -cd $c.gz | /bin/grep $mi;    then options "$ni"
# else options '-t devpts devpts'
# fi

options -t devpts -o newinstance -o ptmxmode=0666 devpts
