# trunk/grimoire/init-scripts.d/init.d/reboot
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 6
# Required-Start: finale
# Should-Start: reboot-kexec
# Short-Description: reboots the box
### END INIT INFO

. /lib/lsb/init-functions

if ! rootfs; then trap - EXIT; exit 0; fi

only start
deny control

start(){
 log_success_msg Rebooting
 /sbin/reboot -d -f -i
}
