# trunk/grimoire/init-scripts.d/init.d/root-rw-retry
#!/bin/bash                                                   

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: root-rw
# Required-Stop: $local_fs                                                      
# Default-Start: S
# Short-Description: provides a second chance for mounting / read-write
### END INIT INFO

if   [ -d /home ] ||
     grep -q boot=IR /proc/cmdline ||
     grep -v rootfs  /proc/mounts  | grep -q " /media/root/[^/]* .*rw"
then exit
fi

. /lib/lsb/init-functions

only start
deny control

start(){ /etc/init.d/root-rw start; }
