# trunk/grimoire/init-scripts.d/init.d/rootname
#!/bin/bash

# Copyright Kyle Sallee 2011, 2012, 2013
# All rights reserved.
# For use with the distribution Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: 1 2 3 4 5
# Short-Description: sets the root file system hostname from /etc/hostname
### END INIT INFO

. /lib/lsb/init-functions

only start configure

configure(){
 if rootfs; then
  log_failure_msg $"This is rootfs."
  log_failure_msg $"rootnames can only be set within a root file system."
  return 1
 fi

 hostname_dialog(){
  BACKTITLE=$"rootname configuration menu"
      TITLE=$"rootname input"
       HELP=$"\
The hostname for assigned to a root file system\n
is also called the rootname.\n
Each root file system's hostname must be unique.\n
The rootname will be saved to file /etc/hostname\n
Please enter the rootname."

  /bin/dialog \
   --backtitle "$BACKTITLE"  \
   --title     "$TITLE"      \
   --timeout 60              \
   --stdout --no-cancel      \
   --inputbox "$HELP" 0 0
 }

 confirm_it(){ /bin/dialog --nocancel --yesno "Missing dot.  Use $NAME as is?" 0 0; }

 local NAME=$( hostname_dialog )
 if   echo "$NAME" | /bin/grep -q "\." || confirm_it
 then echo "$NAME" > /etc/hostname; /bin/rm -f /etc/rootname
 fi
}

start(){
 local REPLY
 log_warning_msg "rootname setting"
 if   [ -f /etc/rootname ]; then read REPLY < /etc/rootname
 elif [ -f /etc/hostname ]; then read REPLY < /etc/hostname
 fi

 if   [ -z "$REPLY" ]
 then REPLY="unknown$RANDOM$RANDOM"
      log_failure_msg "/etc/init.d/rootname configure"
      log_failure_msg "Executing the above command is recommended"
 fi

 echo "$REPLY" > /proc/sys/kernel/hostname 
 log_success_msg "$rootname is now $REPLY"
}
