# trunk/grimoire/init-scripts.d/init.d/run_lock
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: run
# Default-Start: S 1 2 3 4 5
# Short-Description: creates /run/lock
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /run/lock

start(){ [ -d /run/lock ] || mkdir -m 755 /run/lock; }
