# trunk/grimoire/init-scripts.d/init.d/scsi_wait_scan
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: modprobe-coldplug udevsettle
# Default-Start: S
# Short-Description: scsi_wait_scan loads the module
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
deny control
