# trunk/grimoire/init-scripts.d/init.d/shm
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dev
# Default-Start: S 1 2 3 4 5
# Short-Description: mounts tmpfs on /dev/shm
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /dev/shm
options '-o defaults,size=256m,nr_inodes=16k -t tmpfs tmpfs '
