# trunk/grimoire/init-scripts.d/init.d/sorcery
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts tmpfs on /aux/run/tmp
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /aux/run/tmp
options '-o defaults,size=256m,nr_inodes=256k,mode=700 -t tmpfs tmpfs'

stop(){
 active(){
  pgrep -u root '^moil$'
  pgrep -u root '^oust$'
  pgrep -u root '^post$'
  pgrep -u root '^reap$'
  pgrep -u root '^save$'
 }

 busy(){ active | grep -q .; }

 abort_and_wait(){
  touch /aux/run/tmp/halt
  augur queue | sed 's:moil/active/::' | augur abort
  log_warning_msg 'waiting for completion of sorcery'; echo
  log_warning_msg 'forcing power off will cause problems please be patient'
  while busy; do sleep 1; done
  log_success_msg 'sorcery run completed'
 }

 wait_sorcery(){
  if   log_warning_msg 'sorcery run check starting'; ! busy
  then log_success_msg 'sorcery not running'
  else log_warning_msg 'sorcery awaiting completion'; abort_and_wait
  fi
 }

 umount_sorcery(){
  if   log_warning_msg "$NAME unmounting"; umount $NAME
  then log_success_msg "$NAME unmounted"
  else log_success_msg "$NAME not unmounted"
  fi
 }

 wait_sorcery
 umount_sorcery
}
