# trunk/grimoire/init-scripts.d/init.d/sorcery
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts tmpfs on /root/.sorcery
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /root/.sorcery
options '-o defaults,size=256m,nr_inodes=256k,mode=700 -t tmpfs tmpfs'
