# trunk/grimoire/init-scripts.d/init.d/tmp
#!/bin/bash

### BEGIN INIT INFO
# Provides: tmp
# Required-Start: $local_fs
# Should-Start:
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: tmp mounts a 2G tmpfs on /tmp and binds /tmp/.sorcery to /usr/src/sorcery
# Description: tmp mounts a 2G tmpfs on /tmp and binds /tmp/.sorcery to /usr/src/sorcery
### END INIT INFO

. /lib/lsb/init-functions


start()  {
  go()  {
    mount -o defaults,size=2g,nr_inodes=1m,mode=1777 -t tmpfs tmpfs /tmp  &&
    mkdir -p -m 0700 /tmp/.sorcery  /usr/src/sorcery  &&
    mount --bind     /tmp/.sorcery  /usr/src/sorcery
  }

  if    ! sed   's/#.*//'   /etc/fstab  |
          grep  -q          /tmp;  then
    if    log_warning_msg  '/tmp mounting'; go
    then  log_success_msg  '/tmp mounted'
    else  log_failure_msg  '/tmp not mounted'
    fi
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
