# trunk/grimoire/init-scripts.d/init.d/udevd-stop
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: S
# Short-Description: stops udevd running on rootfs
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 # udevadm settle
 /etc/init.d/udevd stop &>/dev/null
 rm -fr /run/udev
}
