# trunk/grimoire/init-scripts.d/init.d/update-mime-database
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: kexec $local_fs udev
# Should-Stop: root-ro
# Should-Stop: tmp
# Default-Stop: 0 6
# Short-Description: updates mime database in mime directories
### END INIT INFO

[ -f /usr/bin/update-mime-database ] || exit 0

. /lib/lsb/init-functions

only stop
deny control

check(){
 local DIR
 while read DIR; do
  log_warning_msg "$DIR check"
    if ! [ -d $DIR ]; then :
  elif ! [ -f $DIR/mime.cache ]; then echo "$DIR"
  elif [ -n "$( /bin/find $DIR -mindepth 2 -newer $DIR/mime.cache )" ]
  then echo "$DIR"
  fi
  log_success_msg "$DIR check"
 done
}

run(){
 local ret=0
 while read DIR; do
  if   log_warning_msg "update-mime-database $DIR execute"; /usr/bin/update-mime-database $DIR &>/dev/null
  then log_success_msg "update-mime-database $DIR execute"
  else log_failure_msg "update-mime-database $DIR execute"; ret=1
  fi
 done
 return $ret
}

stop(){
 echo "/usr/share/mime
/opt/32bit/usr/share/mime" |
 check | run
}
