# trunk/grimoire/init-scripts.d/init.d/usb
#!/bin/bash

### BEGIN INIT INFO
# Provides: usb
# Required-Start: $local_fs
# Should-Start: proc udev
# Default-Start: 1 2 3 4 5
# Short-Description: usb mounts usbfs on /proc/bus/usb
### END INIT INFO

. /lib/lsb/init-functions

NAME=/proc/bus/usb
OPTS='-t usbfs usbfs'

parse "$@"
