# trunk/grimoire/init-scripts.d/init.d/usb
#!/bin/bash

### BEGIN INIT INFO
# Provides: usb
# Required-Start: $local_fs
# Should-Start: proc udev
# Required-Stop: 
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: usb mounts usbfs on /proc/bus/usb
# Description: usb mounts usbfs on /proc/bus/usb
### END INIT INFO

. /lib/lsb/init-functions

SERV='mount'
NAME='/proc/bus/usb'
OPTS='-t usbfs usbfs'

parse "$@"
