# trunk/grimoire/init-scripts.d/init.d/usb
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Required-Start: proc
# Short-Description: mounts usbfs on /proc/bus/usb
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
deny control
name /proc/bus/usb
options '-t usbfs usbfs'
