# trunk/grimoire/init-scripts.d/init.d/usr_src_sorcery_bind
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs
# Required-Stop: proc swap
# Should-Start: tmp
# Should-Stop: tmp
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: binds /tmp/.sorcery to /usr/src/sorcery
### END INIT INFO

case $1 in
start) [[ -d /tmp/.sorcery ]] || mkdir -pm 700 /tmp/.sorcery ;;
esac

. /lib/lsb/init-functions

NAME=/usr/src/sorcery
OPTS='--bind /tmp/.sorcery'
parse "$@"
