# trunk/grimoire/init-scripts.d/init.d/utmp
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: var_lock
# Default-Start: S 1 2 3 4 5
# Short-Description: creates login record /var/run/utmp
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control
name /var/run/utmp

start(){
 run(){
  echo  -n  >    $NAME &&
  chown root:tty $NAME &&
  chmod 664      $NAME
 }

 if   log_warning_msg  "$NAME creating"; run
 then log_success_msg  "$NAME created"
 else log_failure_msg  "$NAME not created"
 fi
}
