# trunk/grimoire/init-scripts.d/init.d/var_run
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: run
# Default-Start: S
# Short-Description: creates bind from /run to /var/run
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control
name /var/run
options '--bind /run'
