# trunk/grimoire/init-scripts.d/init.d/var_lock
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: run_lock
# Default-Start: S
# Short-Description: creates bind from /run/lock to /var/lock
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control
name /var/lock
options '--bind /run/lock'
