# trunk/grimoire/init-scripts.d/init.d/var_lock
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: run $local_fs
# Default-Start: S 1 2 3 4 5
# Short-Description: creates bind from /run/lock to /var/lock
### END INIT INFO

. /lib/lsb/init-functions

only start
deny control
name /var/lock
options '--bind /run/lock'
