# trunk/grimoire/init-scripts.d/init.d/var_tmp
#!/bin/bash

# Copyright Kyle Sallee 2013
# All rights reserved.
# For use with the distribution Modern Magic only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: tmp
# Default-Start: S 1 2 3 4 5
# Short-Description: validates {,/var}/tmp
### END INIT INFO

(( UID == 0 )) || exit 0

. /lib/lsb/init-functions

only start
name /var/tmp

agu(){ [ "1777 0 0" == "$( /bin/stat -c "%a %g %u" "$1" )" ]; }

init_tmp(){
 if   ! [ -d $1 ]; then /bin/mkdir -pm 1777  "$1"
 elif ! agu  $1  ; then /bin/chmod     1777  "$1"
                        /bin/chown root:root "$1"
 fi
}

start(){
 local r

 if   log_warning_msg "/var/tmp/ validating"; init_tmp /var/tmp/
 then log_success_msg "/var/tmp/ validated"
 else log_failure_msg "/var/tmp/ validation failed"; r=1
 fi

 if   log_warning_msg "/tmp/ validating"; init_tmp /tmp/
 then log_success_msg "/tmp/ validated"
 else log_failure_msg "/tmp/ validation failed"; r=1
 fi

 return $r
}
