# trunk/grimoire/openntpd.d/init.d/openntpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Provides: ntpd
# Required-Start: $network
# Required-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: openntpd synchronizes system time with a server
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/ntpd
NAME=ntpd

start(){
 if   log_warning_msg "$NAME starting"; start_daemon -b $SERV -s
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

status(){
 pgrep ^$NAME >/dev/null  &&
 log_warning_msg "$NAME is running"  ||
 log_warning_msg "$NAME is stopped"
}

stop(){
 log_warning_msg "Stopping $NAME"
 pkill ^$NAME >/dev/null 2>&1     &&
 log_success_msg "$NAME stopped"  ||
 log_failure_msg "$NAME not stopped"
}

parse "$@"
