# trunk/grimoire/module-init-tools.d/init.d.on/modprobe
#!/bin/bash

### BEGIN INIT INFO
# Provides: modprobe
# Required-Start: sys
# Should-Start:
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: modprobe loads modules during boot
# Description: modprobe loads modules specified in /etc/modules during boot
### END INIT INFO

. /lib/lsb/init-functions

mp(){
 if   [ -f          /etc/modules ]
 then sed 's:#.*::
	   /^ *$/d' /etc/modules /sys/bus/pci/devices/*/modalias
 else cat                        /sys/bus/pci/devices/*/modalias
 fi |
 while read
 do modprobe $REPLY
 done
}

start(){
 log_warning_msg 'Modules loading'; mp
 log_success_msg 'Modules loaded'
}

help(){ echo "Usage: $0 {start}" 1>&2;}

case $1 in
  start) start ;;
try-restart) : ;;
      *) help  ;;
esac
