# trunk/grimoire/gizmod.d/init.d/gizmod
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $syslog run
# Required-Stop: kill
# Should-Stop: $syslog
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: gizmod provides the gizmod server
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
pid_file /var/run/gizmod.pid
server /usr/bin/gizmod

start(){
 log_warning_msg 'gizmod starting'
 $SERV &
 sleep 0.1
 echo $! > $PIDF
 log_success_msg 'gizmod started'
}
