# trunk/grimoire/net-tools.d/hostname
#!/bin/sh

### BEGIN INIT INFO
# Provides: hostname
# Required-Start: 
# Required-Stop:  
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: hostname sets the box's name and starts loopback
# Description: hostname sets a box's name and starts the loopback device
### END INIT INFO

/bin/hostname  -F  /etc/hostname
/sbin/ifconfig  lo  127.0.0.1
