# trunk/grimoire/net-tools.d/network
#!/bin/sh

### BEGIN INIT INFO
# Provides: $network
# Required-Start: hostname $local_fs
# Required-Stop:  $local_fs
# Default-Start: 3 4 5
# Default-Stop:  0 1 2 6
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

# Edit and uncomment the appropriate lines

DEVICE=
IP=
BROADCAST=
NETMASK=
GATEWAY=
MODULE=

case $1 in
    start)  
            log_warning_msg "networking about to start"
            echo      "$1ing networking."
            # modprobe  $MODULE
            # ifconfig  $DEVICE $IP broadcast $BROADCAST netmask $NETMASK
            # route     add default gateway $GATEWAY
            # dhcpcd    $DEVICE
            log_success_msg "networking started"
            ;;
              
     stop)  log_warning_msg "networking about to end"
            # ifconfig  $DEVICE down
            # modprobe  -r  $MODULE
            log_success_msg "networking is ended"
            ;;

  restart)  $0  stop  &&
            $0  start
            ;;

        *)  echo  "Usage: $0 {start|stop|restart}"
            ;;
esac
