# trunk/grimoire/util-linux-ng.d/init.d/hwclock
#!/bin/bash

### BEGIN INIT INFO
# Provides: localtime
# Required-Start: $local_fs
# Should-Start: $local_fs
# Required-Stop:  $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: Sets and saves the system clock by the hardware clock.
# Description: Sets and saves the system clock by the hardware clock.
### END INIT INFO

. /lib/lsb/init-functions

start()  {  /sbin/hwclock  --hctosys  --utc;  }
stop()   {  /sbin/hwclock  --systohc  --utc;  }
help()   {  echo "Usage: $0 {start}" 1>&2;    }

case  $1  in
        start)  start  ;;
         stop)  stop   ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
