# trunk/grimoire/NVIDIA-Linux-x86_64.d/init.d.on/agpgart
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: agpgart loads the module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
