# trunk/grimoire/NVIDIA-Linux-x86_64.d/init.d.on/nvidiactl
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs nvidiafb agpgart nvidia
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: nvidiactl creates device nodes.
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 run(){
  for ((CX=0;CX<8;CX++)); do
  mknod            /dev/nvidia$CX c 195 $CX; done
  mknod            /dev/nvidiactl c 195 255
  chmod 660        /dev/nvidia*
  chown root:video /dev/nvidia*
 }

 log_warning_msg 'nvidiactl starting'; run 2>/dev/null
 log_warning_msg 'nvidiactl started'
}
