# trunk/grimoire/dhcp.d/init.d/dhcpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs $network run
# Required-Stop:  $local_fs $network
# Default-Start: S
# Default-Stop:  0 6
# Short-Description: provides dynamic hosting control protocol daemon
### END INIT INFO

# set IF to the interface device in
# /etc/init.d/conf.d/dhcpd
# of the root file system
# used to generate the rootfs initramfs images.

. /lib/lsb/init-functions

fast
only start stop status
server /sbin/dhcpd
config /etc/dhcpd.conf
options "-q ${IF:=eth0}"

# pid_file /var/run/dhcpd.pid

start(){
 /sbin/ip route add 255.255.255.255 dev $IF
 /bin/mkdir -pm 755 /var/db/
 /bin/rm -f /var/db/dhcpd.leases
 /bin/touch /var/db/dhcpd.leases
}

stop(){
 /sbin/ip route del 255.255.255.255 dev $IF
 /bin/rm -f /var/db/dhcpd.leases
}
