# trunk/grimoire/lm_sensors.d/init.d.on/sensors
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Should-Start: sensors_modules
# Required-Stop: syslog kill 
# Default-Start: 2 3 4 5
# Short-Description: sensors sets hardware limits specified in /etc/sensors.conf
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
server /usr/bin/sensors
config /etc/sensors3.conf
options -s
