# trunk/grimoire/lm_sensors.d/init.d.on/sensord
#!/bin/bash

### BEGIN INIT INFO
# Provides: sensord
# Required-Start: sensors run
# Should-Start: 
# Required-Stop: syslog kill 
# Should-Stop: 
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: sensord is lm_sensors hardware health monitoring daemon
# Description: sensord starts the lm_sensors hardware health monitoring
# daemon which log the health status into syslog or a round-robin database. 
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/sensord.pid
SERV=/usr/bin/sensord
NAME=sensord

# Configuration of sensord

SCAN=1m		# interval between scanning for sensor alarms
LAG=30m		# interval between logging  all sensor readings

OPTS="-i $SCAN -l $LAG"

parse "$@"
