# trunk/grimoire/lm_sensors.d/init.d.on/sensors
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Should-Start: sensors_modules
# Required-Stop: syslog kill 
# Default-Start: 2 3 4 5
# Short-Description: sensors sets hardware limits specified in /etc/sensors.conf
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/bin/sensors
CONF=/etc/sensors3.conf
OPTS=-s

help(){ echo "Usage: $0 {start}" 1>&2; }

parse "$@"
