# trunk/grimoire/fcron.d/init.d.on/fcron
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Provides: $cron
# Required-Start: $local_fs run
# Required-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: fcron is a vixie cron replacement
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/fcron.pid
SERV=/usr/sbin/fcron
CONF=/etc/fcron.conf

start(){
 export TERM=dumb
 if   ! [ -f $CONF ]; then
      log_warning_msg "fcron not started; $CONF does not exist"
 elif log_warning_msg "fcron starting"; start_daemon -p $PIDF $SERV
 then log_success_msg "fcron started"
 else log_failure_msg "fcron not started"; status
 fi
}

parse "$@"
