# trunk/grimoire/courier-authlib.d/init.d.on/authdaemond
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: $local_fs
# Required-Stop: $local_fs
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: courier authdaemon
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/authdaemond
CONF=/etc/authlib/authdaemonrc

restart(){ OPTS=stop; stop && sleep 1; OPTS=start; start; }

help(){ echo "Usage: $0 {restart|force-reload|start|stop}" 1>&2; }

case "$@" in
 start) OPTS=start ;;
  stop) OPTS=stop ;;
esac

parse "$@"
