# trunk/grimoire/systemd.d/init.d/ptmx
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Should-Start: udevsettle
# Default-Start: S 1 2 3 4 5
# Short-Description: creates /dev/ptmx symbolic link if it disappears
### END INIT INFO

. /lib/lsb/init-functions

only start
fast

start(){
 if ! [   -h          /dev/ptmx ]
 then  rm -f          /dev/ptmx
       ln -s pts/ptmx /dev/ptmx
 fi
}
