# trunk/grimoire/systemd.d/init.d/udevtrigger
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: udevd
# Default-Start: S 1 2 3 4 5
# Short-Description: udevtrigger turns on dynamic /dev population and module loading
# Description: udevtrigger requests kernel device uevents.
# Coldplug events are used by udevd to popualte /dev
### END INIT INFO

. /lib/lsb/init-functions

# to view subsystems
# cat /lib/udev/rules.d/* | tr ' ' '\n' | grep SUBSYSTEM | sort -u

only start
cgroup none
server /sbin/udevadm
if   [ -d /home ]
then options trigger
else options trigger --subsystem-nomatch=bluetooth --subsystem-nomatch=drm --subsystem-nomatch=graphics --subsystem-nomatch=power_supply --subsystem-nomatch=printer --subsystem-nomatch=sound --subsystem-nomatch=video4linux
fi
