# trunk/grimoire/systemd.d/init.d/udevd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: udev
# Required-Start: dev-cp sys
# Should-Start: cgroup cgroup-init-rootfs cgroup-init-root
# Required-Stop: $local_fs kill
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: udevd provides dynamic /dev populates and module loading
### END INIT INFO

# [ -d /lib/modules/ ] || exit 0

. /lib/lsb/init-functions

only start stop restart try-restart status
server /lib/systemd/systemd-udevd
options --daemon
