# trunk/grimoire/systemd.d/init.d/udevtrigger
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: udevd
# Default-Start: S 1 2 3 4 5
# Short-Description: udevtrigger turns on dynamic /dev population and module loading
# Description: udevtrigger requests kernel device uevents.
# Coldplug events are used by udevd to populate /dev
### END INIT INFO

[ -d /lib/modules/ ] || exit 0

. /lib/lsb/init-functions

only start
cgroup none
server /bin/udevadm
options trigger
