# trunk/grimoire/wpa_supplicant.d/init.d/var_run_wpa_priv
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Default-Start: S
# Short-Description: creates directory /var/run/wpa_priv
### END INIT INFO

. /lib/lsb/init-functions

only start
name /var/run/wpa_priv

start(){
 run(){ mkdir -pm 7505 $NAME && chown root:wpapriv $NAME; }

 if   [ -d $NAME ]; then return; fi
 if   log_warning_msg "$NAME creating"; run
 then log_success_msg "$NAME created"
 else log_failure_msg "$NAME not created"; return 1
 fi
}
