# trunk/grimoire/wpa_supplicant.d/init.d/wpa_priv
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: wpa_priv
# Required-Start: proc var_run_wpa_priv
# Should-Start: udev
# Required-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Starts wpa_priv for all wireless interfaces
### END INIT INFO

. /lib/lsb/init-functions

pnw(){
 [ -f          /proc/net/wireless ] &&
 tr ' ' '\n' < /proc/net/wireless |
 sed 's/://p;d' |
 sed 's/^/wext:/'
}

only start stop restart status
name wpa_priv
server /sbin/wpa_priv
options "-B -c /var/run/wpa_priv -P /var/run/wpa_priv.pid $( pnw )"
