# trunk/grimoire/immune-security-tomoyo.d/init.d/tomoyo-transition-sorcerer
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: tomoyo-domain
# Default-Start: S 1 2 3 4 5
# Short-Description: loads Sorcerer specific policies for Tomoyo
### END INIT INFO

. /lib/lsb/init-functions

only start try-restart

start(){
 copy(){ cp -a /etc/tomoyo/sorcerer/transition/* /domain/; }

 if   log_warning_msg "tomoyo domain transition programs copying"; copy
 then log_success_msg "tomoyo domain transition programs copied"
 else log_failure_msg "tomoyo domain transition programs not copied"; return 1
 fi
}
