# trunk/grimoire/immune-security-tomoyo.d/init.d/tomoyo-auditd
#!/bin/bash

# Copyright Kyle Sallee 2011, 2012, 2013
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: securityfs
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: start the tomoyo logger
### END INIT INFO

. /lib/lsb/init-functions

only start status stop
server /sbin/tomoyo-auditd

if [ "$1" == start ]; then
 [ -d /var/log/        ] || /bin/mkdir -pm 755 /var/log/
 [ -d /var/log/tomoyo/ ] || /bin/mkdir -pm 700 /var/log/tomoyo/
fi
