# trunk/grimoire/immune-security-tomoyo.d/init.d/tomoyo-auditd
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: proc
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: start the tomoyo logger
### END INIT INFO

. /lib/lsb/init-functions

only start restart status stop try-restart
server /sbin/tomoyo-auditd

if [ "$1" == start ]; then
 [ -d           /var/log/tomoyo ] ||
 mkdir -pvm 700 /var/log/tomoyo
fi
