# trunk/grimoire/immune-security-tomoyo.d/init.d/tomoyo-domain
#!/bin/bash

# Copyright Kyle Sallee 2011, 2012, 2013, 2014
# All rights reserved.

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts a 256K tmpfs on /domain
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /domain
options "-o defaults,size=256k,nr_inodes=256,mode=755 -t tmpfs tmpfs"
