# trunk/grimoire/kdebase-workspace.d/init.d/kdm
#!/bin/bash

### BEGIN INIT INFO
# Provides: kdm gui_login
# Required-Start: $local_fs run ICE-unix
# Should-Start: $network $remote_fs udev xfs
# Required-Stop: kill
# Should-Stop: $remote_fs $network initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: kdm provides an x11 login.
# Description: kdm provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/kdm.pid
SERV=/usr/bin/kdm
NAME=kdm

parse "$@"
