# trunk/grimoire/bind.d/init.d/named
#!/bin/bash

### BEGIN INIT INFO
# Provides: $named
# Required-Start: $network
# Required-Stop: $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: named provides name to IP address translation
# Description: named provides name to IP address translation
### END INIT INFO

. /lib/lsb/init-functions

JAIL=/var/chroot/bind
PIDF=$JAIL/var/run/named.pid
SERV=/usr/sbin/named
NAME=named

status()  {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}

start()  {
  log_warning_msg  "$NAME about to start"
  modprobe capability 2>/dev/null
  if    start_daemon     -p $PIDF $SERV -u bind -t $JAIL
  then  log_success_msg  "$NAME started"
  else  log_failure_msg  "$NAME not started"
        status
  fi
}

stop()  {
  log_warning_msg "$NAME about to stop"
  rndc stop
  log_success_msg "$NAME stopped"
#  if    killproc        -p $PIDF -s 15  $SERV
#  then  log_success_msg "$NAME stopped"
#  else  log_failure_msg "$NAME not stopped"
#        status
#        return 0
#  fi
}

restart() {  stop  &&  sleep  1;  start;  }

try_restart() {
  status  &&  restart
  return 0
}

reload()  {
  rndc reload
#  PIDS="$( pidofproc -p $PIDF $SERV )"
#  if    [ -n            "$PIDS" ]
#  then  kill  -s SIGHUP  $PIDS  &&
#        log_success_msg "$NAME configuration reloaded"
#  fi
}

help()  {
  echo "Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}

case  $1  in
       restart)      restart  ;;
   try-restart)  try_restart  ;;
        reload)  reload       ;;
  force-reload)  reload       ;;
         start)  start        ;;
        status)  status       ;;
          stop)  stop         ;;
             *)  help         ;;
esac
