# trunk/grimoire/sorcery.d/init.d.on/immune
#!/bin/bash

### BEGIN INIT INFO
# Provides: immune
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: kill 
# Should-Stop: $remote_fs
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: immune provides a computer immune system
# Description: immune is capable of detecting and fixing critical
#	library breaks, neglect, reporting near full partitions and
#	compressing files in the background
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/immune.pid
SERV=/usr/sbin/immune
NAME=immune

parse
