# trunk/grimoire/net-tools.d/init.d/network
#!/bin/sh

### BEGIN INIT INFO
# Provides: $network
# Required-Start: hostname $local_fs
# Required-Stop:  $local_fs
# Default-Start: 3 4 5
# Default-Stop:  0 1 2 6
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

if    [  -x  /etc/init.d/networking.sh  ]
then  log_warning_msg  "Please edit /etc/init.d/network and remove /etc/init.d/networking.sh"
      exec   /etc/init.d/networking.sh  "$@"
fi

. /lib/lsb/init-functions

# Edit and uncomment the appropriate lines

DEVICE=
IP=
BROADCAST=
NETMASK=
GATEWAY=
MODULE=

case $1 in
    start)  
            log_warning_msg "networking about to start"
            # modprobe  $MODULE
            # ifconfig  $DEVICE $IP broadcast $BROADCAST netmask $NETMASK
            # route     add default gateway $GATEWAY
            # dhcpcd    $DEVICE
            log_success_msg "networking started"
            ;;
              
     stop)  log_warning_msg "networking about to stop"
            # ifconfig  $DEVICE down
            # modprobe  -r  $MODULE
            log_success_msg "networking is stopped"
            ;;

  restart)  $0  stop  &&
            $0  start
            ;;

        *)  echo  "Usage: $0 {start|stop|restart}"
            ;;
esac
