# trunk/grimoire/net-tools.d/init.d/network
#!/bin/bash

### BEGIN INIT INFO
# Provides: $network
# Required-Start: $local_fs hostname
# Required-Stop:  $local_fs
# Default-Start: 3 4 5
# Default-Stop:  0 1 2 6
# Short-Description: network configures the box for static IP address networking
# Description: network configures the box for static IP address networking
### END INIT INFO

. /lib/lsb/init-functions

# Edit and uncomment the appropriate lines

DEVICE=
IP=
BROADCAST=
NETMASK=
GATEWAY=
MODULE=

case $1 in
    start)  if    [  -n  "$BROADCAST"  ]  &&
                  [  -n  "$DEVICE"     ]  &&
                  [  -n  "$GATEWAY"    ]  &&
                  [  -n  "$IP"         ]  &&
                  [  -n  "$NETMASK"    ]
            then
              log_warning_msg "networking about to start"
              if    [  -n                  "$MODULE"  ]  &&
                    !  lsmod  |  grep -qw  "$MODULE" 
              then  modprobe               "$MODULE"
              fi
              ifconfig  $DEVICE $IP broadcast $BROADCAST netmask $NETMASK
              route     add default gateway $GATEWAY
              log_success_msg "networking started"
            fi
            ;;
              
     stop)  log_warning_msg "networking about to stop"
            [  -n     "$DEVICE"  ]  &&
            ifconfig   $DEVICE  down
            log_success_msg "networking has stopped"
            ;;

  restart)  $0  stop  &&
            $0  start
            ;;

        *)  echo  "Usage: $0 {start|stop|restart}"
            ;;
esac
