# trunk/grimoire/cryptsetup.d/init.d.on/dm-crypt
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dm-mod
# Should-Start: udev
# Required-Stop: $syslog
# Default-Start: S
# Short-Description: dm-crypt loads the dm-crypt module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
