# trunk/grimoire/cryptsetup.d/init.d.on/luks
#!/bin/bash

### BEGIN INIT INFO
# Provides: luks
# Required-Start: dmsetup
# Should-Start: udev vgchange
# Required-Stop: halt reboot
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: luks finds luks encrypted devices and maps them to /dev/mapper/
### END INIT INFO

. /lib/lsb/init-functions

NAME=luks

start(){
 run(){
  local DEV
  for   DEV in $( sed 's/#.*//;/noauto/d;;s/\t/\n/g;s/ /\n/g' /etc/fstab |
                  sed 's,/dev/mapper/,/dev/,p;d' ); do
   cryptsetup isLuks   /dev/$DEV &&
   cryptsetup luksOpen /dev/$DEV $DEV
  done
 }

 log_warning_msg "$NAME checking for luks devices"; run <> /dev/console
 log_success_msg "$NAME luks mappings created"
}


stop(){
 log_warning_msg "$NAME closing luks devices"
 find /dev/mapper -type b |
 while read
 do cryptsetup luksClose ${REPLY##*/}
 done
 log_success_msg "$NAME luks devices closed."
}

try_restart(){ :; }

parse "$@"
