# trunk/grimoire/util-linux.d/init.d/devices
#!/bin/bash

### BEGIN INIT INFO
# Provides: devices
# Required-Start: 
# Should-Start: 
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: devices mounts devfs on /devices for old boxes
# Description: devices mounts devfs on /devices for old boxes
### END INIT INFO

. /lib/lsb/init-functions


start()  {
  log_warning_msg  'devfs on /devices mounting'
  if    [  -d  /devices  ]  &&
        mount  -n  -t devfs  devfs /devices
  then  log_success_msg  'devfs on /devices mounted'
  else  log_failure_msg  'devfs on /devices not mounted'
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
