# trunk/grimoire/util-linux.d/init.d/dmesg
#!/bin/bash

### BEGIN INIT INFO
# Provides: dmesg
# Required-Start: $local_fs $syslog
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: dmesg prints kernel ring buffer to /var/log/dmesg
# Description: dmesg prints kernel ring buffer to /var/log/dmesg
### END INIT INFO

start() { /bin/dmesg -c > /var/log/dmesg; }
help()  { echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
 start)  start ;;
     *)  help  ;;
esac
