# trunk/grimoire/util-linux.d/init.d/localtime
#!/bin/bash

### BEGIN INIT INFO
# Provides: swap
# Required-Start: $local_fs
# Required-Stop:  $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: localtime setup the localtime
# Description: localtime
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg "localtime about to activate";  
	#/sbin/hwclock  --hctosys  --localtime
	log_warning_msg "configure /etc/init.d/localtime please!"
	false
  log_status_msg  $? "localtime active"
}

stop()  {
  log_warning_msg  "swap about to deactivate";
  log_success_msg  "swap deactivated"
}

help()  {
  echo "Usage: $0 {start|stop}" 1>&2
}

case  $1  in
         start)  start        ;;
          stop)  stop         ;;
             *)  help         ;;
esac
