# trunk/grimoire/util-linux.d/init.d/localtime
#!/bin/bash

### BEGIN INIT INFO
# Provides: localtime
# Required-Start: $local_fs
# Required-Stop:  $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: Invoke this at boot when hardware clock stores localtime instead of UTC.
# Description: localtime sets the system time according to the
#	hardware clock and assumes the hardware clock contains localtime
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  /sbin/hwclock  --hctosys  --localtime
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
         start)  start        ;;
             *)  help         ;;
esac
