# trunk/grimoire/util-linux.d/init.d/swap
#!/bin/bash

### BEGIN INIT INFO
# Provides: swap
# Required-Start: $local_fs
# Required-Stop:  $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: swap turns on swaps specified in /etc/fstab
# Description: swapon and swapoff
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg  "swap about to activate";  swapon  -a
  log_success_msg  "swap active"
}

stop()  {
  log_warning_msg  "swap about to deactivate";  swapoff -a
  log_success_msg  "swap deactivated"
}

help()  {  echo "Usage: $0 {start|stop}" 1>&2;  }

case  $1  in
        start)  start  ;;
         stop)  stop   ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
