# trunk/grimoire/util-linux.d/init.d/utmp-wtmp
#!/bin/bash

### BEGIN INIT INFO
# Provides: utmp-wtmp
# Required-Start: $local_fs $syslog
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: utmp-wtmp create login records /var/run/utmp /var/log/wtmp
# Description:       utmp-wtmp create login records /var/run/utmp /var/log/wtmp
### END INIT INFO

start()  {
    /usr/bin/touch            /var/run/utmp  /var/log/wtmp  2>/dev/null
  ( /bin/chown      root:tty  /var/run/utmp  /var/log/wtmp  2>/dev/null
    /bin/chmod      0664      /var/run/utmp  /var/log/wtmp  2>/dev/null ) &
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
         start)  start        ;;
             *)  help         ;;
esac
