# branches/grimoire/grimoire/netkit-combo/xinetd.d/fingerd
service finger
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= nobody
	group		= tty
	server		= /usr/sbin/in.fingerd
	server_args	= 
	only_from	= 127.0.0.1
	per_source	= 16
	log_type	= FILE /var/log/xinetd/fingerd
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
