# trunk/grimoire/tor.d/init.d/tor
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs veth
# Should-Start: hwclock ntpd
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: tor is an onion router daemon
## With a default configuration with relay or bridge enabled
## then remote computers connect to port 9001
## Therefore, the network firewall
## must be configured to permit connection
## to this computer on port 9001 TCP.
## This computer's firewall must be adjusted by editing either
## /+/base/etc/port.forward or
## /etc/port.receive in the root file system running tor
## must be edited to allow inbound connections to the VETH's port 9001
## If running only as a client rather than a bridge or relay
## then no additional adjustments are required beyond editing
## /etc/tor/torrc
## Unless the computer is protected by a fascist firewall.
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
config /etc/tor/torrc
domain limited
server /usr/bin/tor
