# trunk/grimoire/dbus.d/init.d/var_run_dbus
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Default-Start: S 1 2 3 4 5
# Short-Description: creates directory /var/run/dbus
### END INIT INFO

. /lib/lsb/init-functions

only start
name /var/run/dbus

start(){
 run(){ mkdir -pm 775 $NAME && chown messagebus:messagebus $NAME; }

 if   [ -d $NAME ]; then return; fi
 if   log_warning_msg "create"; run
 then log_success_msg "create"
 else log_failure_msg "create"; return 1
 fi
}
