# trunk/grimoire/thttpd.d/thttpd
#!/bin/sh

# chkconfig: 345 90 10
# description: Starts and stops the thttpd web server.

.  /lib/lsb/init-functions

case $1 in
    start)  echo  "Starting thttpd web server."
            /usr/sbin/thttpd  -C  /etc/thttpd.conf
            ;;
     stop)  echo  "Stopping thttpd web server."
            killproc  thttpd
            ;;
  restart)  echo  "Restarting thttpd web server."
            killproc  thttpd
            /usr/sbin/thttpd
            ;;
esac
