# trunk/grimoire/device-mapper.d/init.d.on/dm-mod
#!/bin/bash

### BEGIN INIT INFO
# Provides: dm-mod
# Required-Start:
# Should-Start: udev
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: dm-mod loads the device mapper module
# Description: dm-mod loads the device mapper module
### END INIT INFO

. /lib/lsb/init-functions

NAME=dm-mod

parse "$@"
