# trunk/grimoire/device-mapper.d/init.d.on/dmsetup
#!/bin/bash

### BEGIN INIT INFO
# Provides: dmsetup
# Required-Start:
# Should-Start: udev
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: dmsetup creates the /dev/mapper/control device
# Description: dmsetup creates the /dev/mapper/control device
### END INIT INFO

. /lib/lsb/init-functions

NAME=dmsetup
OPTS='mknodes'
SERV=/sbin/dmsetup

start(){
 if   log_warning_msg "$NAME $OPTS"; $SERV $OPTS
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
       start) start ;;
 try-restart) : ;;
           *) help  ;;
esac
