# trunk/grimoire/device-mapper.d/init.d.on/dmsetup
#!/bin/bash

### BEGIN INIT INFO
# Provides: dmsetup
# Required-Start: dm-mod
# Should-Start: udev
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: dmsetup creates the /dev/mapper/control device
# Description: dmsetup creates the /dev/mapper/control device
### END INIT INFO

. /lib/lsb/init-functions

NAME=dmsetup
OPTS='mknodes'
SERV=/sbin/dmsetup

parse "$@"
