# trunk/grimoire/kdenetwork-v3.d/init.d/lisa
#!/bin/bash

### BEGIN INIT INFO
# Provides: lisa
# Required-Start: $network run
# Should-Start: $remote_fs udev xfs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: lisa provides the LAN Information Server
# Description: lisa provides the LAN Information Server
#	LISa is intended to provide a kind of "network neighbourhood"
#	but only relying on the TCP/IP protocol stack, no smb or whatever.
### END INIT INFO

. /lib/lsb/init-functions

NAME=lisa
PIDF=/var/run/lisa.pid
SERV=/opt/kde3/bin/lisa
OPTS='-q'

parse "$@"
