# trunk/grimoire/alsa-utils.d/init.d.on/alsactl
#!/bin/bash

### BEGIN INIT INFO
# Provides: alsactl
# Should-Start: snd-mixer-oss snd-pcm-oss snd-seq-oss
# Default-Start: 1 2 3 4 5
# Short-Description: alsactl restores the audio mixer volume settings
### END INIT INFO

. /lib/lsb/init-functions

NAME=alsactl
SERV=/usr/sbin/alsactl
OPTS=restore

parse "$@"
