# trunk/grimoire/alsa-utils.d/init.d.on/alsactl
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Should-Start: snd-mixer-oss snd-pcm-oss snd-seq-oss
# Default-Start: 1 2 3 4 5
# Short-Description: alsactl restores the audio mixer volume settings
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/alsactl
OPTS=restore
CONF=/etc/asound.state

parse "$@"
