# trunk/grimoire/alsa-utils.d/init.d.on/snd-seq-oss
#!/bin/bash

### BEGIN INIT INFO
# Provides: snd-seq-oss
# Should-Start: snd-mixer-oss
# Default-Start: 1 2 3 4 5
# Short-Description: snd-seq-oss loads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME=snd-seq-oss

parse "$@"
