# trunk/grimoire/alsa-utils.d/init.d.on/snd-pcm-oss
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Should-Start: modprobe udev
# Default-Start: 1 2 3 4 5
# Short-Description: snd-pcm-oss loads the module
### END INIT INFO

. /lib/lsb/init-functions

parse "$@"
