# trunk/grimoire/mailutils/xinetd.d/imap4d-stunnel
service imaps
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/sbin/stunnel
	server_args		= -l /usr/sbin/imap4d
        # server		= /usr/sbin/imapd
	per_source		= 4
	log_type		= FILE /var/log/xinetd/imap4d-stunnel
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
