# trunk/grimoire/mailutils/xinetd.d/pop3d-stunnel
service pop3s
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/sbin/stunnel
	server_args		= -l /usr/sbin/pop3d
	# server		= /usr/sbin/pop3d
	per_source		= 4
	log_type		= FILE /var/log/xinetd/pop3d-stunnel
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
