# trunk/grimoire/suspend-utils.d/init.d.on/resume
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: udev
# Default-Start: S
# Short-Description: resumes from suspend to RAM or disk.
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 if   log_warning_msg "resuming"; echo | /sbin/resume
 then log_success_msg "resumed"
 else log_failure_msg "resume failed"
 fi
}
