# trunk/grimoire/courier-imap.d/init.d/courier-pop3
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: pop3
# Required-Start: $local_fs $network authdaemon
# Required-Stop: $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: courier-imap, POP3 daemon
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
server /usr/libexec/pop3d.rc
name courier-pop3

start(){
 log_warning_msg "$NAME starting"; $SERV start
 log_success_msg "$NAME started"
}

stop(){
 log_warning_msg "$NAME stopping"; $SERV stop
 log_success_msg "$NAME stopped"
}

cd /usr/libexec
