# trunk/grimoire/VirtualBox.d/init.d.on/vboxdrv
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs $syslog
# Should-Start: udev
# Default-Start: 2 3 4 5
# Short-Description: loads the vboxdrv module
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
