# trunk/grimoire/VirtualBox.d/init.d.on/vboxnetflt
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs $syslog vboxdrv
# Should-Start: udev
# Default-Start: 2 3 4 5
# Short-Description: loads the vboxnetflt module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
