# trunk/grimoire/rsyslog.d/init.d.on/rsyslogd
#!/bin/bash

### BEGIN INIT INFO
# Provides: rsyslogd $syslog
# Required-Start: hostname $local_fs run
# Should-Start: hwclock
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: rsyslogd is a replacement for sysklog and syslog
# Description: rsyslogd is a near replacement for sysklog, sysklogd and syslog
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/rsyslogd.pid
SERV=/usr/sbin/rsyslogd
NAME="rsyslogd"

parse "$@"
