# trunk/grimoire/coreutils.d/init.d.on/hostname
#!/bin/bash

### BEGIN INIT INFO
# Provides: hostname
# Required-Start: proc
# Required-Stop: 
# Default-Start: S
# Default-Stop: 
# Short-Description: hostname sets the box's name
# Description: hostname sets the box's name
### END INIT INFO

. /lib/lsb/init-functions

start() {
 log_warning_msg "hostname setting"
 /bin/hostname -F /etc/hostname
 log_success_msg "hostname set"
}

case $1 in
 start) start ;;
esac
