# trunk/grimoire/lighttpd.d/init.d/lighttpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: lighttpd is a web server
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/lighttpd.pid
SERV=/usr/sbin/lighttpd
OPTS='-f /etc/lighttpd.conf'

parse "$@"
