# trunk/grimoire/lighttpd.d/init.d/lighttpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: lighttpd is a web server
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status

pid_file    /var/run/lighttpd.pid
server      /usr/sbin/lighttpd
config      /etc/lighttpd/lighttpd.conf
options '-f /etc/lighttpd/lighttpd.conf'
