# trunk/grimoire/ltmodem.d/init.d.on/ltserial
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs ltmodem
# Should-Start: udev
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Short-Description: ltserial loads the module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
