# trunk/grimoire/postfix.d/init.d.on/newaliases
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: set-permissions
# Default-Start: 2 3 4 5
# Short-Description: updates the postfix's alias database
### END INIT INFO

. /lib/lsb/init-functions

only start
server /usr/bin/newaliases
name newaliases
