# trunk/grimoire/postfix.d/init.d.on/set-permissions
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: var_spool_postfix
# Default-Start: 2 3 4 5
# Short-Description: creates proper directories inside /var/spool/postfix
### END INIT INFO

. /lib/lsb/init-functions

only start
server /usr/sbin/postfix
options set-permissions
