# trunk/grimoire/rsyslog.d/init.d/rsyslogd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: $syslog
# Required-Start: hostname $local_fs run
# Should-Start: hwclock ntpd
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: rsyslogd is a replacement for sysklog and syslog
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
#pid_file /var/run/rsyslogd.pid
server /usr/sbin/rsyslogd
