# trunk/grimoire/cpufrequtils.d/init.d/cpufreq_gouvernor
#!/bin/dash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: S
# Short-Description: cpufreq_gouvernor sets the box's cpu freq gouvernor
### END INIT INFO

. /lib/lsb/init-functions

start(){
 log_warning_msg "cpufreq_gouvernor setting"
 if [ -f           /etc/cpufreq_gouvernor ]; then
   read   REPLY < /etc/cpufreq_gouvernor
   cpufreq-set -g ${REPLY}
   log_success_msg "cpufreq_gouvernor is now $REPLY"
 else
   log_failure_msg "/etc/cpufreq_gouvernor don't  exist"
 fi
}

parse "$@"
