# trunk/grimoire/glibc.d/init.d.on/var_run_nscd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: run
# Default-Start: 2 3 4 5
# Short-Description: var_run_nscd mounts tmpfs on /var/run/nscd
### END INIT INFO

. /lib/lsb/init-functions

name /var/run/nscd
config /etc/nscd.conf
options '-o defaults,size=16k,nr_inodes=4,mode=700 -t tmpfs tmpfs'
