# trunk/grimoire/glibc.d/init.d.on/var_db_nscd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: 2 3 4 5
# Short-Description: var_db_nscd mounts tmpfs on /var/db/nscd
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /var/db/nscd
config /etc/nscd.conf
options '-o defaults,size=256m,nr_inodes=64,mode=700 -t tmpfs tmpfs'
