# trunk/grimoire/glibc.d/init.d.on/var_run_nscd
#!/bin/bash

### BEGIN INIT INFO
# Provides: var_run_nscd
# Required-Start: run
# Should-Start:
# Required-Stop: 
# Default-Start: 2 3 4 5
# Default-Stop:
# Short-Description: var_run_nscd mounts tmpfs on /var/run/nscd
### END INIT INFO

. /lib/lsb/init-functions

NAME=/var/run/nscd
CONF=/etc/nscd.conf
OPTS='-o defaults,size=16k,nr_inodes=4,mode=700 -t tmpfs tmpfs'

parse "$@"
