# trunk/grimoire/install.d/init.d.on/lvm
#!/bin/bash

### BEGIN INIT INFO
# Provides: lvm
# Required-Start: tmp
# Should-Start: find_ir $local_fs udev
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: lvm copies /etc/lvm to /tmp then binds /tmp/lvm to /etc/lvm
# Only pvcreate a vgchange from LVM2 work properly if /etc/lvm is read only.
### END INIT INFO

. /lib/lsb/init-functions

start(){
 on_ir(){ grep -q boot=IR /proc/cmdline; }

 run(){
  cp -a /etc/lvm /tmp
  mount --bind /tmp/lvm /etc/lvm
 }

 if   on_ir &&
      log_warning_msg '/etc/lvm bind starting'; run
 then log_success_msg '/etc/lvm bind completed'
 fi
}


help(){ echo "Usage: $0 {start}" 1>&2; }


case $1 in
       start) start ;;
 try-restart) :     ;;
           *) help  ;;
esac
