# trunk/grimoire/install.d/init.d.on/welcome
#!/bin/bash

### BEGIN INIT INFO
# Provides: welcome
# Required-Start:
# Should-Start: dmesg hostname utmp-wtmp
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: welcome instructs the SA how to begin installation
# Description:       welcome instructs the SA how to begin installation
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  welcome()  {
    echo  $"Welcome to Sorcerer."
    echo  $"Please login as root."
    echo  $"Press enter for password."
    echo  $"/install"
    echo  $"to begin a menu driven installation."
  }

  log_warning_msg  "welcome starting";  welcome  >  /dev/console
  log_success_msg  "welcome completed"
}


help()  {  echo "Usage: $0 {start}" 1>&2;  }


case  $1  in
        start)  start ;;
  try-restart)  :     ;;
            *)  help  ;;
esac
