# trunk/grimoire/oss.d/init.d.on/ossdevlinks
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: osssound
# Default-Start: 1 2 3 4 5
# Short-Description: ossdevlinks creates legacy symbolic links for OSS device nodes
### END INIT INFO

. /lib/lsb/init-functions

name ossdevlinks
server /usr/sbin/ossdevlinks

start(){
 if   log_warning_msg "$SERV run"; $SERV
 then log_success_msg "$SERV exited successfully"
 else log_failure_msg "$SERV exited with an error"
 fi
}

stop(){ :; }
