# trunk/grimoire/oss.d/init.d.on/ossdetect
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: ossdetect detects soundcards and writes config to /etc/oss/installed_drivers
### END INIT INFO

if [ -f /etc/oss/installed_drivers ]; then exit 0; fi

. /lib/lsb/init-functions

only start reload
name ossdetect
server /usr/sbin/ossdetect
options -i

start(){
 if   log_warning_msg "About to run $SERV to detect audio devices."; $SERV $OPTS
 then log_success_msg "$SERV exited successfully"
 else log_failure_msg "$SERV exited with an error"
 fi
}

reload(){ rm -f /etc/oss/installed_drivers; start; }
