# trunk/grimoire/oss.d/init.d.on/ossmixer
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: osssound
# Default-Start: 1 2 3 4 5
# Short-Description: ossmixer loads previously saved mixer settings
### END INIT INFO

. /lib/lsb/init-functions

only start
name savemixer
server /usr/sbin/savemixer
options -L

start(){
 if   log_warning_msg "$SERV $OPTS running"; $SERV $OPTS
 then log_success_msg "$SERV exited successfully"
 else log_failure_msg "$SERV exited with an error"
 fi
}
