# trunk/grimoire/oss.d/init.d.on/osssound
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: ossdetect
# Default-Start: 1 2 3 4 5
# Short-Description: soundon loads OSS modules
### END INIT INFO

. /lib/lsb/init-functions

deny try-restart

start(){
 NAME=soundon
 SERV=/usr/sbin/soundon
 if   log_warning_msg "About to run $SERV to load OSS modules"; $SERV
 then log_success_msg "$SERV exited successfully"
 else log_failure_msg "$SERV exited with an error"
 fi
}

stop(){
 NAME=soundoff
 SERV=/usr/sbin/soundoff
 if   log_warning_msg "About to unload OSS modules";  $SERV
 then log_success_msg "$SERV excited successfully"
 else log_failure_msg "$SERV exited with an error"
 fi
}
