# trunk/grimoire/gdm.d/init.d.on/gdm
#!/bin/bash

### BEGIN INIT INFO
# Provides: gui_login
# Required-Start: $network dbus run ICE-unix
# Should-Start: $remote_fs udev xfs dbus
# Required-Stop: kill
# Should-Stop: $remote_fs initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: gdm provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

#PIDF=/var/run/gdm.pid
SERV=/usr/sbin/gdm
CONF=/etc/gdm/custom.conf
LOGF=/var/log/gdm.log

sidofproc(){
 ps axo sess,comm | grep gdm-binary | cut -d ' ' -f1 | head -n 1
}

status(){
 local SIDS="$( sidofproc )"
 if ["$PIDS" ]
 then log_warning_msg "$NAME is running"; return 0
 else log_warning_msg "$NAME is not running"; return 3
 fi
}


stop(){
 local SIDS=$( sidofproc )
 if   log_warning_msg "$NAME stopping"; pkill -15 -s $SIDS
 then log_success_msg "$NAME stopped"
 else log_failure_msg "$NAME not stopped"; status
 fi
}

start(){
 if   log_warning_msg "$NAME starting"; start_daemon -b -l $LOGF $SERV $OPTS
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

parse "$@"
