# trunk/grimoire/hotplug
   VERSION=( "2004_09_23" "2004_03_29" "2004_01_05" "2003_08_05" )
    IGNORE="^S2000"
  CATEGORY="utility"
 ATTRIBUTE="console"
    SOURCE="pub/linux/utils/kernel/hotplug/hotplug-$VERSION.tar.bz2"
       URL="$KERNEL_URL"
       REQ="init-functions"
  HOMEPAGE="http://linux-hotplug.sourceforge.net"
  ESTIMATE="5"
      DESC="hotplug lets you plug in new devices and use them immediately."

build() {
   TH="touch /var/lock/subsys/hotplug"
MKDIR="s:$TH:mkdir -p /var/lock/subsys; $TH:"

sed  -i  "$MKDIR
2i### BEGIN INIT INFO
2i# Provides: hotplug
2i# Required-Start: \$local_fs
2i# Required-Stop:  kill
2i# Default-Start: 1 2 3 4 5
2i# Default-Stop: 0 6
2i# Short-Description: hotplug loads modules for new devices 
2i# Description: hotplug creates devices nodes and loads 
2i#	kernel drivers for new devices immediately
2i### END INIT INFO"  \
etc/rc.d/init.d/hotplug

prepare_install  &&
make    install
}

post_install() {
  /lib/lsb/install_initd /etc/init.d/hotplug
  true
}

current() { grep  -q  "kill"  /etc/init.d/hotplug; }
