# trunk/grimoire/multipath-tools.d/init.d.on/multipath-flush
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: multipathd
# Should-Start: udev 
# Default-Start: S
# Short-Description: multipath-flush flushes all unused multipath device maps
### END INIT INFO

. /lib/lsb/init-functions

NAME=multipath
SERV=/sbin/multipath
OPTS=-F

parse "$@"
