# trunk/grimoire/multipath-tools.d/init.d.on/multipathd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: mapper
# Should-Start: udev 
# Default-Start: S
# Short-Description: multipathd allows a host to access a LU by multiple paths
# Description: mulitpathd allows a host to access 
#	Logical Units (LU) by multiple paths. The most
#	common multipathed environment today is Fibre
#	Channel (FC) Storage Area Network (SAN).
### END INIT INFO

. /lib/lsb/init-functions

name	multipathd
server	/sbin/multipathd
parse "$@"
