# trunk/grimoire/udev.d/init.d/udev
#!/bin/bash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: 
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg  "udev about to activate"
  /bin/mount  -n  -t proc   proc   /proc
  /bin/mount  -f  -t proc   proc   /proc
  /bin/mount      -t sysfs  sysfs  /sys 
  /bin/mount      -t ramfs  ramfs  /dev
  [           -d         /lib/udev/devices  ]  &&
  /bin/cp     -a         /lib/udev/devices/*  /dev
  /sbin/udevstart
  /sbin/udevd --daemon
  log_success_msg  "udev started"
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
  start)  start  ;;
      *)  help   ;;
esac
