# trunk/grimoire/udev.d/init.d/udev
#!/bin/bash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: 
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg  "udev about to activate"
  mount  -n  -t proc   none  /proc
  mount      -t sysfs  none  /sys 
  /etc/init.d/start_udev
	udevd --daemon
  #pgrep  udevd >/dev/null 2>&1  ||  /sbin/udevsend
  log_success_msg  "udev started"
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
         start)  start        ;;
             *)  help         ;;
esac
