# trunk/grimoire/udev.d/init.d/udevd-restart
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: udev
# Default-Start: 1 2 3 4 5
# Short-Description: restarts udevd on the real root file system
### END INIT INFO

. /lib/lsb/init-functions

only start
fast

start(){
 run (){
  udevadm settle
  /etc/init.d/udevd restart
  udevadm trigger
  udevadm settle
  [ -h /dev/ptmx ] ||
  /etc/init.d/ptmx start
 }

 run &>/dev/null &
 disown
}
