# trunk/grimoire/udev.d/init.d/udev
#!/bin/bash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: 
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  log_warning_msg  "udev about to activate"
  /bin/mount  -n  -t proc   none  /proc
  /bin/mount  -f  -t proc   none  /proc
  /bin/mount      -t sysfs  none  /sys 
  /bin/mount      -t ramfs  none  /dev
  /bin/cp     -a  /lib/udev/devices/*  /dev
  /sbin/udevstart
  /sbin/udevd --daemon
  log_success_msg  "udev started"
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
         start)  start        ;;
             *)  help         ;;
esac
