# trunk/grimoire/udev.d/init.d/udevsettle
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: udev
# Required-Start: udevd udevtrigger
# Default-Start: S
# Short-Description: udevsettle waits for queued udev events to be handled
### END INIT INFO

. /lib/lsb/init-functions

only start configure status
cgroup none
server /sbin/udevadm
options "settle --timeout ${TIMEOUT:-60}"

status(){
 log_warning_msg "udevsettle timeout set for ${TIMEOUT:-60} seconds";
# if   [ -n "$( pidofproc /sbin/udevadm )" ]
# then log_warning_msg "udevsettle is running"
# else log_warning_msg "udevsettle is not running"
# fi
}

configure(){

timeout_dialog(){
 BACKTITLE=$"udevsettle Configuration Menu"
     TITLE=$"Timeout Selection"
      HELP=$"udevsettle terminates after udevd completes queued events.
udevsettle requires less than 10 seconds.
Running longer hints at the potential of a problem.
The timeout can be set for less than 60 seconds
when udev processing hangs or lingers
and yet causes no adverse effect upon booting.
If lowering the timeout causes file systems
to not fsck or to not mount
then the timeout was set too low.

Most SAs never change the udevsettle timeout value.
Change this only when one already possesses the knowledge
required for fixing the problem if the box fails to boot.

Please select the amount of seconds for timeout duration."

 dialog \
  --backtitle "$BACKTITLE"  \
  --title     "$TITLE"      \
  --timeout 60              \
  --stdout --no-cancel      \
  --default-item 60         \
  --menu   "$HELP" 0 0 0    \
  300 seconds 240 seconds 180 seconds 90 seconds 60 seconds 30 seconds 15 seconds
}

 TIMEOUT=$( timeout_dialog )
 if   [ -n "$TIMEOUT" ]
 then echo "TIMEOUT=$TIMEOUT" > /etc/init.d/conf.d/udevsettle
 fi
}
