# trunk/grimoire/cryptmount.d/init.d.on/cryptmount
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs
# Should-Start: udev
# Required-Stop: halt reboot
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: cryptmount configures encrypted devices at boot.
### END INIT INFO

. /lib/lsb/init-functions

name cryptmount
server /bin/cryptmount
config /etc/default/cryptmount

CM_BOOTDV=
CM_BOOTSW=
CM_BOOTFS=

status(){
 local ON OFF
 for TARGET in $CM_BOOTDV $CM_BOOTSW
 do [ -b "/dev/mapper/${TARGET}" ] && ON+="$TARGET " || OFF+="$TARGET "
 done
 if [ "$OFF" ]
 then log_warning_msg "$OFF not active"; return 3;
 else log_warning_msg "$ON active";
 fi
}

start(){
 run(){
  if ! status 2>/dev/null; then
   if [ "$CM_BOOTDV" ]; then $SERV --prepare $CM_BOOTDV; fi &&
   if [ "$CM_BOOTSW" ]; then $SERV --swapon  $CM_BOOTSW; fi &&
   if [ "$CM_BOOTFS" ]; then $SERV --mount   $CM_BOOTFS; fi
  fi
 }

 if   log_warning_msg "$NAME starting"; run
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

stop(){
 run(){
  if [ "$CM_BOOTFS" ]; then $SERV --unmount $CM_BOOTFS; fi
  if [ "$CM_BOOTSW" ]; then $SERV --swapoff $CM_BOOTSW; fi
  if [ "$CM_BOOTDV" ]; then $SERV --release $CM_BOOTDV; fi
  $SERV -n
 }

 if   log_warning_msg "$NAME stopping"; $SERV $FPTS -n
 then log_success_msg "$NAME stopped"
 else log_failure_msg "$NAME not stopped"; status
 fi
}

if [ -f $CONF ]; then . $CONF; fi
