# trunk/grimoire/cryptmount.d/init.d.on/dm-crypt
#!/bin/bash

### BEGIN INIT INFO
# Provides: dm-crypt
# Required-Start:
# Should-Start: udev
# Required-Stop: $syslog
# Default-Start: S
# Default-Stop:
# Short-Description: dm-crypt loads the dm-crypt module
# Description: dm-crypt loads the dm-crypt module
### END INIT INFO

. /lib/lsb/init-functions

NAME=dm-crypt

parse "$@"
