# trunk/grimoire/fuse.d/init.d.on/fuse
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Should-Start: udev
# Required-Stop: halt reboot
# Default-Start: S
# Short-Description: fuse loads the module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
