# trunk/grimoire/fuse.d/init.d.on/fusectl
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: fuse
# Should-Start: udev
# Required-Stop: halt reboot
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: fusectl mounts fusectl on /sys/fs/fuse/connections
### END INIT INFO

. /lib/lsb/init-functions

NAME=/sys/fs/fuse/connections
OPTS='-t fusectl fusectl'

parse "$@"
