# trunk/grimoire/tomoyo-tools.d/init.d/tomoyo-domain
#!/bin/bash

# Copyright Kyle Sallee 2011
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S
# Short-Description: mounts a 64M tmpfs on /domain
### END INIT INFO

. /lib/lsb/init-functions

only start restart status
name /domain
options "-o defaults,size=64m,nr_inodes=1m,mode=755 -t tmpfs tmpfs"
