# trunk/grimoire/hal.d/init.d.on/hald
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: dbus run $syslog
# Required-Stop: kill
# Should-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: hald provides present system devices through dbus system bus
### END INIT INFO


. /lib/lsb/init-functions

PIDF=/var/run/haldaemon.pid
SERV=/usr/sbin/hald
OPTS='--daemon=yes --use-syslog'

parse "$@"
