# trunk/grimoire/hal.d/init.d.on/hald
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: dbus run $syslog var_cache_hald
# Required-Stop: kill
# Should-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: hald provides present system devices through dbus system bus
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
pid_file /var/run/haldaemon.pid
server /usr/sbin/hald
options '--daemon=yes --use-syslog'
