# trunk/grimoire/kbd.d/init.d/kbd
#!/bin/bash

### BEGIN INIT INFO
# Provides: kbd
# Required-Start: $local_fs $syslog
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: kbd loads the keyboard map and console font
# Description:       kbd loads the keyboard map and console font
### END INIT INFO


start()  {
  /bin/loadkeys  us
  /usr/bin/setfont  default8x16
}


help()  {  echo "Usage: $0 {start}" 1>&2;  }


case  $1  in
        start)  start  ;;
         stop)  :      ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
