# trunk/grimoire/radvd.d/init.d/radvd-ipv6-port-forward
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: proc $network
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: enables port forwarding of IPv6 packets for all network interfaces
### END INIT INFO

. /lib/lsb/init-functions

start(){
 log_warning_message "echo 1 > /proc/sys/net/ipv6/conf/all/forwarding"
 echo 1 > /proc/sys/net/ipv6/conf/all/forwarding
 log_success_message "echo 1 > /proc/sys/net/ipv6/conf/all/forwarding complete"
}

stop(){ :;}

parse "$@"
