# trunk/grimoire/immune-boot-sentient.d/init.d/boot-sentient
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Short-Description: Transitions from Sorcerer to Sorcerer Sentient upon boot
### END INIT INFO

if [ "$1" == try-restart ]]; then exit 5; fi
if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 5; fi

. /lib/lsb/init-functions

server /usr/sbin/boot-sentient
