# trunk/grimoire/ttyd.d/init.d/ttyd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: provides dynamic login daemon with root selection
### END INIT INFO

. /lib/lsb/init-functions

only start stop status
server /sbin/ttyd
domain user
