# trunk/grimoire/samba.d/init.d/nmbd
#!/bin/bash

### BEGIN INIT INFO
# Provides: nmbd
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: nmbd provides a NetBIOS name server
# Description: nmbd provides a NetBIOS name server
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/nmbd.pid
SERV=/usr/sbin/nmbd
NAME=nmbd
OPTS='-D'

parse
