# trunk/grimoire/samba.d/init.d/nmbd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $network run
# Should-Start: $remote_fs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: nmbd provides a NetBIOS name server
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
server /usr/sbin/nmbd
config /etc/samba/smb.conf
options -D
#pid_file /var/run/nmbd.pid
