# trunk/grimoire/msyslog.d/init.d.on/msyslog
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Provides: $syslog
# Required-Start: hostname $local_fs msyslog-logs run
# Should-Start: hwclock
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: msyslog is a replacement for sysklog and syslog
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/syslogd.pid
SERV=/usr/sbin/syslogd
CONF=/etc/syslog.conf
OPTS='-i linux -i unix'

parse "$@"
