# trunk/grimoire/iproute2.d/init.d.on/hostname
#!/bin/bash

### BEGIN INIT INFO
# Provides: hostname
# Required-Start: $local_fs
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: hostname sets the box's name and starts loopback
# Description: hostname sets a box's name and starts the loopback device
### END INIT INFO

start() {
 read < /etc/hostname
 echo $REPLY > /proc/sys/kernel/hostname
 /sbin/ip addr replace 127.0.0.1 dev lo
 /sbin/ip link set lo up
}

case $1 in
 start) start ;;
esac
