# trunk/grimoire/gdk-pixbuf.d/init.d/gdk-pixbuf-query-loaders
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 510
# Required-Stop: $local_fs
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: updates giomodule.cache
### END INIT INFO

# Copyright 2011 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

. /lib/lsb/init-functions

only stop try-restart force-reload status

try_restart(){
/etc/init.d/gdk-pixbuf-query-loaders-32 stop
/etc/init.d/gdk-pixbuf-query-loaders-64 stop
}

stop(){
/etc/init.d/gdk-pixbuf-query-loaders-32 stop
/etc/init.d/gdk-pixbuf-query-loaders-64 stop
}

status(){
/etc/init.d/gdk-pixbuf-query-loaders-32 status &&
/etc/init.d/gdk-pixbuf-query-loaders-64 status
}

force_reload(){
/etc/init.d/gdk-pixbuf-query-loaders-32 force-reload
/etc/init.d/gdk-pixbuf-query-loaders-64 force-reload
}
