# trunk/grimoire/initramfs.d/init.d/initramfs-xen
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: kexec $local_fs root-ro udev
# Should-Stop: tmp
# Default-Stop:
# Short-Description: creates initramfs that xen can load
### END INIT INFO

# Copyright 2007 through 2013 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

[ -f /usr/sbin/xl ] || exit 0

. /lib/lsb/init-functions

only stop reload force-reload status
name /boot/xen/initrf
INIT=initramfs-bin

allpaths(){ /bin/sed -nr "p;:a;s:(.+)/.*:\1:p;ta" | /bin/sort -u; }

exists(){ while read; do [[ -e $REPLY ]] && echo "$REPLY"; done; }

every(){
 local lver="$(
  /bin/find /boot -mindepth 2 -maxdepth 2 -type f -name linux -printf '%T@\t%h\n' |
  /bin/sort -gr |
  /bin/sed '\,/boot/rescue,d; s:.*/::' |
  /bin/head -n 1
 )"

 cd /
 /bin/find /boot/$lver/modules /boot/initrf/ -type f |
 while read
 do    /bin/gzip -cd "$REPLY" | /bin/cut -b2-1024 | /bin/cpio -t --quiet
 done
}

initrf_update(){
 local               TEMPIF=/tmp/initramfs.$$.$RANDOM$RANDOM
 /bin/rm    -rf     $TEMPIF
 /bin/mkdir -pm 700 $TEMPIF/boot/xen

 local COMMON="$TEMPIF/$NAME"

 every | allpaths | exists |
 /bin/cpio -o -H newc --quiet |
 /bin/gzip -1 > $COMMON
 /bin/chmod 600 $COMMON

 /bin/cp -a  $TEMPIF/boot /
 /bin/rm -fr $TEMPIF
 /bin/sync
}

initrf_present(){ [ -f $NAME ] || return 1; }

initrf_current(){
 if   /bin/find /boot/initrf/ -type f -newer $NAME | /bin/grep -q .
 then return 1
 fi

 if [ -f /var/log/init/failed/$INIT/stop ]; then return 1; fi
} 

run(){
 if   log_warning_msg "update"; initrf_update
 then log_success_msg "update"
 else log_failure_msg "update"; false
 fi
}

status(){
 if   log_warning_msg "check"; initrf_present && initrf_current
 then log_success_msg "check status:     current"
 else log_warning_msg "check status: not current"; return 1
 fi
}

stop(){ status || run; }
force_reload(){ run; }
reload(){ DONE=true; stop; }

export LC_ALL=C
