# trunk/grimoire/initramfs.d/init.d/initramfs
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: kexec $local_fs root-ro udev
# Should-Stop: tmp
# Default-Stop: 0 6
# Short-Description: creates boot filesystems
### END INIT INFO

# Copyright 2007 through 2013 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

[ -d /boot/ ] || exit 0

. /lib/lsb/init-functions

only stop reload force-reload status
name initramfs

initrf_update(){
 local r=0

 for REPLY in /etc/init.d/initramfs-*
 do  bash $REPLY $1 || r=1
 done

 return $r
}

run(){
 if   log_warning_msg "update"; initrf_update $1
 then log_success_msg "update"
 else log_failure_msg "update"; false
 fi
}

status(){
 for REPLY in /etc/init.d/initramfs-*
 do  /bin/bash $REPLY status || return 1
 done
}

stop(){ status || run stop; }
force_reload(){ run force-reload; }
reload(){ stop; }

export LC_ALL=C
