# trunk/grimoire/initramfs.d/init.d/initramfs-sbin
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: kexec $local_fs root-ro udev
# Should-Stop: tmp
# Default-Stop:
# Short-Description: creates boot filesystems
### END INIT INFO

# Copyright 2007 through 2013 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

. /lib/lsb/init-functions

only stop reload force-reload status
name /boot/initrf/sbin
INIT=initramfs-sbin

allpaths(){ /bin/sed -nr "p;:a;s:(.+)/.*:\1:p;ta" | /bin/sort -u; }

exists(){ while read; do [[ -e $REPLY ]] && echo "$REPLY"; done; }

every(){
 /bin/find "/${NAME##*/}" |
 /bin/sed "\
\,^/sbin/dhcpd$,d
\,^/sbin/dhcrelay$,d
\,^/sbin/mount\.ecryptfs$,d
\,^/sbin/mount\.smbfs$,d
\,^/sbin/sln$,d"
}

initrf_update(){

 local               TEMPIF=/tmp/initramfs.$$.$RANDOM$RANDOM
 /bin/rm    -rf     $TEMPIF
 /bin/mkdir -pm 700 $TEMPIF/boot/initrf

 local  COMMON="$TEMPIF/$NAME"

 cd /
 every |  allpaths | exists |
 /bin/cut -b2-1024 |
 /bin/cpio -o -H newc --quiet |
 /bin/gzip -1 > $COMMON
 /bin/chmod 600 $COMMON

 [       -f /boot/initrf ] &&
 /bin/rm -f /boot/initrf

 /bin/cp -a  $TEMPIF/boot /
 /bin/rm -fr $TEMPIF
 /bin/sync
}

initrf_present(){ [ -f $NAME ] || return 1; }

initrf_current(){

 if   /bin/find "/${NAME##*/}" -newer $NAME | /bin/grep -q .
 then return 1
 fi

 if [ -f /var/log/init/failed/$INIT/stop ]; then return 1; fi

} 

run(){
 if   log_warning_msg "update"; initrf_update
 then log_success_msg "update"
 else log_failure_msg "update"; false
 fi
}

status(){
 if   log_warning_msg "check"; initrf_present && initrf_current
 then log_success_msg "check status:     current"
 else log_warning_msg "check status: not current"; return 1
 fi
}

stop(){ status || run; }
force_reload(){ run; }
reload(){ DONE=true; stop; }

export LC_ALL=C
