# trunk/grimoire/linux.d/init.d/tune_kernel_io
#!/bin/sh

### BEGIN INIT INFO
# Provides: 
# Required-Start: $local_fs
# Default-Start: S
# Short-Description: tune kernel for particular use
# Description: tune kernel for particular use
### END INIT INFO
. /lib/lsb/init-functions



#see /usr/src/linux/Documentation/filesystems/proc.txt

# Maximum time, in seconds, of hard drive spindown time that you are
# confortable with. Worst case, it's possible that you could lose this
# amount of work if your battery fails you while in laptop mode.
MAX_AGE=600

# Dirty synchronous ratio.  At this percentage of dirty pages the process which
# calls write() does its own writeback
DIRTY_RATIO=40


#
# Allowed dirty background ratio, in percent.  Once DIRTY_RATIO has been
# exceeded, the kernel will wake pdflush which will then reduce the amount
# of dirty memory to dirty_background_ratio.  Set this nice and low, so once
# some writeout has commenced, we do a lot of it.
#
DIRTY_BACKGROUND_RATIO=5


# This must be adjusted manually to the value of HZ in the running kernel
# on 2.4, until the XFS people change their 2.4 external interfaces to work in
# centisecs. This can be automated, but it's a work in progress that still needs
# some fixes. On 2.6 kernels, XFS uses USER_HZ instead of HZ for external
# interfaces, and that is currently always set to 100. So you don't need to
# change this on 2.6.
XFS_HZ=${XFS_HZ:-'100'}
AGE=$((100*$MAX_AGE))
XFS_AGE=$(($XFS_HZ*$MAX_AGE))


#############################################################################

case "$1" in
	start)
	  log_warning_msg "kernel tunning"

		# (For 2.4 and early 2.6.)
		# This only needs to be set, not reset -- it is only used when
		# laptop mode is enabled.
		  #echo $XFS_AGE > /proc/sys/vm/pagebuf/lm_flush_age
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_sync_interval
		# (A couple of early 2.6 laptop mode patches had these.)
		# The same goes for these.
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_age_buffer
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_sync_interval
		# (2.6.6)
		# But not for these -- they are also used in normal
		# operation.
		  #echo $XFS_AGE > /proc/sys/fs/xfs/age_buffer
		  #echo $XFS_AGE > /proc/sys/fs/xfs/sync_interval
		# (2.6.7 upwards)
		# And not for these either. These are in centisecs,
		# not USER_HZ, so we have to use $AGE, not $XFS_AGE.
		  #echo $AGE > /proc/sys/fs/xfs/age_buffer_centisecs
		  #echo $AGE > /proc/sys/fs/xfs/xfssyncd_centisecs
		  #echo 3000 > /proc/sys/fs/xfs/xfsbufd_centisecs

		#2.4 kernels
			#echo 1					> /proc/sys/vm/laptop_mode
			#echo "30 500 0 0 $AGE $AGE 60 20 0"	> /proc/sys/vm/bdflush
		
		#2.6 kernels
		#laptop_mode is the time in second to wait when a write is performed to group
		#several write in one call
			#echo 5					> /proc/sys/vm/laptop_mode
			#echo "$AGE"				> /proc/sys/vm/dirty_writeback_centisecs
			#echo "$AGE"				> /proc/sys/vm/dirty_expire_centisecs
			#echo "$DIRTY_RATIO"			> /proc/sys/vm/dirty_ratio
			#echo "$DIRTY_BACKGROUND_RATIO"		> /proc/sys/vm/dirty_background_ratio

    log_success_msg "tunning finished"
		;;
	stop)
	  log_warning_msg "kernel untunning"
		# kernel default dirty buffer age
		AGE=30
		UPDATE=5
		DIRTY_BACKGROUND_RATIO=10
		DIRTY_RATIO=40
		XFS_AGE_BUFFER=15
		XFS_SYNC_INTERVAL=30
		XFS_BUFD_INTERVAL=1
		
		AGE=$((100*$MAX_AGE))
		XFS_AGE=$(($XFS_HZ*$MAX_AGE))

		# (For 2.4 and early 2.6.)
		# This only needs to be set, not reset -- it is only used when
		# laptop mode is enabled.
		  #echo $XFS_AGE > /proc/sys/vm/pagebuf/lm_flush_age
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_sync_interval
		# (A couple of early 2.6 laptop mode patches had these.)
		# The same goes for these.
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_age_buffer
		  #echo $XFS_AGE > /proc/sys/fs/xfs/lm_sync_interval
		# (2.6.6)
		# But not for these -- they are also used in normal
		# operation.
		  #echo $XFS_AGE > /proc/sys/fs/xfs/age_buffer
		  #echo $XFS_AGE > /proc/sys/fs/xfs/sync_interval
		# (2.6.7 upwards)
		# And not for these either. These are in centisecs,
		# not USER_HZ, so we have to use $AGE, not $XFS_AGE.
		  #echo $AGE > /proc/sys/fs/xfs/age_buffer_centisecs
		  #echo $AGE > /proc/sys/fs/xfs/xfssyncd_centisecs
		  #echo $XFS_BUFD_INTERVAL > /proc/sys/fs/xfs/xfsbufd_centisecs

		#2.4 kernels
			#echo 0					> /proc/sys/vm/laptop_mode
			#echo "30 500 0 0 $AGE $AGE 60 20 0"	> /proc/sys/vm/bdflush
		
		#2.6 kernels
			#echo 0					> /proc/sys/vm/laptop_mode
			#echo "$AGE"				> /proc/sys/vm/dirty_writeback_centisecs
			#echo "$AGE"				> /proc/sys/vm/dirty_expire_centisecs
			#echo "$DIRTY_RATIO"			> /proc/sys/vm/dirty_ratio
			#echo "$DIRTY_BACKGROUND_RATIO"		> /proc/sys/vm/dirty_background_ratio
    log_success_msg "untunning finished"
		;;
	*)
		echo "Usage: $0 {start|stop}" 2>&1
		exit 1
		;;

esac

exit 0
