# trunk/grimoire/evms.d/init.d.on/evms_activate
#!/bin/bash

### BEGIN INIT INFO
# Provides: evms_activate
# Required-Start: dmsetup proc fstab sys
# Should-Start: udev
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: evms_activate activates evms volumes
### END INIT INFO

. /lib/lsb/init-functions

NAME=evms_activate
OPTS='-l /etc/init.d/log.d/evms-engine -d warning'
SERV=/sbin/evms_activate

# A properly installed box already has /dev/evms/device_names in /etc/fstab

evms_fstab(){
 sed -i 's:\t: :g; \, / ,d' /etc/fstab
 find /dev/evms -maxdepth 1 -type b -printf '%P\n' |
 LC_ALL=C sort -r |
 while read
 do sed -i "s:/dev/$REPLY :/dev/evms/$REPLY :" /etc/fstab
 done
}

start(){
 if   log_warning_msg "$NAME $OPTS"; $SERV $OPTS # && evms_fstab
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
       start) start ;;
     restart) : ;;
 try-restart) : ;;
           *) help  ;;
esac
