# trunk/grimoire/NVIDIA-Linux-x86.d/init.d.on/nouveau
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: nouveau unloads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME='-r nouveau'

parse "$@"
