# trunk/grimoire/NVIDIA-Linux-x86.d/init.d.on/noveau
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: noveau unloads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME='-r noveau'

parse "$@"
