# trunk/grimoire/NVIDIA-Linux-x86.d/init.d.on/nvidia
#!/bin/bash

### BEGIN INIT INFO
# Provides: sshd
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 2 3 4 5
# Default-Stop: 
# Short-Description: nvidia loads module and creates device nodes.
# Description: nvidia loads module and creates device nodes.
### END INIT INFO

. /lib/lsb/init-functions

NAME=nvidia

start(){
 run(){
  rmmod    nvidiafb
  modprobe agpgart
  modprobe nvidia
  for ((CX=0;CX<8;CX++)); do
  mknod            /dev/nvidia$CX c 195 $CX; done
  mknod            /dev/nvidiactl c 195 255
  chmod 0660       /dev/nvidia*
  chown root:video /dev/nvidia*
 }

 log_warning_msg "$NAME starting"; run 2>/dev/null
 log_warning_msg "$NAME started"
}

help(){ echo"Usage: $0 {start}" 1>&2; }

case $1 in
 start) start ;;
 try-restart) : ;;
 *) help ;;
esac
