# trunk/grimoire/immune-sentinel-retry.d/init.d.on/sentinel-retry
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: retries previously failed compilation
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

DELAY=30	# minutes
TRIES=1

. /lib/lsb/init-functions

fast
only start stop restart try-restart status configure
server /usr/sbin/sentinel-retry
options DELAY TRIES

configure(){
 input(){ dialog --stdout --no-cancel --timeout 60  --title $"$NAME configuration menu" --inputbox "$1" 10 50 $2; }

 DELAY=$( input $"Minutes before retrying failed compile?" ${DELAY:=1} )
 TRIES=$( input $"Retry attempts?" ${TRIES:=1} )

 if   (( DELAY > 0 )) && (( TRIES > 0 ))
 then echo "TRIES=$TRIES; DELAY=$DELAY" > /etc/init.d/conf.d/$NAME
 fi
 try_restart
}
