# trunk/grimoire/openssh.d/init.d.on/ssh-keygen
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network empty run
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Short-Description: ssh-keygen generates missing host keys
### END INIT INFO

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

server		/sbin/ssh-keygen
config		/etc/ssh/sshd_config

start(){
 check_key(){
  if   [ -f $1 ] && [ -f $1.pub ]; then return; fi
  if   log_warning_msg "$NAME generating $1"; ssh-keygen -t $2 -f $1 -N ""
  then log_success_msg "$NAME generated  $1"
  else log_failure_msg "$NAME failed     $1"
  fi
 }

 if ! [ -f $CONF ]; then
  log_warning_msg "$NAME requires $CONF"
  return
 fi

 check_key /etc/ssh/ssh_host_key     rsa1
 check_key /etc/ssh/ssh_host_rsa_key rsa
 check_key /etc/ssh/ssh_host_dsa_key dsa
 check_key /etc/ssh/ssh_host_ecdsa_key  ecdsa
}

help(){ echo "Usage: $0 {start}" 1>&2; }
