# trunk/grimoire/shishi.d/init.d/shishid
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $network run
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: shishid serves requests for tickets.
### END INIT INFO

. /lib/lsb/init-functions

pid_file /var/run/shishid.pid
server /usr/sbin/shishid
name shishid
options '-u shishi'

start(){
 run(){ start_daemon -b -p "$PIDF" $SERV $OPTS &>/var/log/shishid; }
 if   log_warning_msg "$NAME starting"; run
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

reload(){
 PIDS="$( pidofproc -p $PIDF $SERV )"
 if   [ -n          "$PIDS" ]
 then kill -s SIGHUP $PIDS &&
      log_success_msg "$NAME configuration reloaded"
 fi
}
