# trunk/grimoire/sorcerer-iso9660.d/init.d/untar-usr
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Short-Description: Extracts /usr/usr.tar.xz
### END INIT INFO

. /lib/lsb/init-functions

if   ! rootfs
then echo "$0 can only be executed on the rootfs"; trap - EXIT; exit 1
fi

only start
deny control

start(){
 local txz=/usr/usr.tar.xz

 if    [ -f $txz ]; then
  if   log_warning_msg "$txz extracting"; tar -Pxf $txz && rm -f $txz /etc/init.d/untar-usr
  then log_success_msg "$txz extracted"
  else log_failure_msg "$txz extraction failed"; return 1
  fi
 fi
}
