# trunk/grimoire/leafnode/xinetd.d/leafnode
service nntp
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= news
	only_from		= localhost
	server			= /usr/sbin/leafnode
        per_source              = 8
	server_args		=
        log_type                = FILE /var/log/xinetd/leafnode
        log_on_success          = PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
