# trunk/grimoire/dhcpcd.d/init.d.on/var_lib_dhcpcd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: var_lib_dhcpcd mounts tmpfs on /var/lib/dhcpcd
### END INIT INFO

. /lib/lsb/init-functions

name /var/lib/dhcpcd
options '-o defaults,size=64k,nr_inodes=8,mode=0755 -t tmpfs tmpfs'
