# trunk/grimoire/postfix.d/init.d/postfix
#!/bin/bash

### BEGIN INIT INFO
# Provides: postfix $mail
# Required-Start: $network
# Should-Start: $remote_fs
# Required-Stop: $network
# Should-Stop: $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: postfix is a MTA
# Description: postfix is a mail transport agent
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/postfix
NAME=postfix
PROCNAME=/usr/libexec/postfix/master
OPTS=start

status(){
 PIDS="$( pidofproc $PROCNAME )"
 if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
 elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                      else log_warning_msg "$NAME is not running"; return 3
 fi
}

start(){
 if   log_warning_msg "$NAME starting"; newaliases && start_daemon $SERV $OPTS
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

stop(){
 log_warning_msg "$NAME stopping"; $SERV stop
 log_success_msg "$NAME stopped"
}

parse "$@"
