# trunk/grimoire/immune-auto-network.d/init.d/auto-network
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: hostname loopback run var_lock_run
# Should-Start: $networks modprobe udev udevsettle
# Required-Stop:  $local_fs
# Should-Stop: kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: dhcpcd configuration of wired and wireless interfaces
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
server /sbin/auto-network
