# tags/start/grimoire/netkit-combo/xinetd.d/ftpd
service ftp
{
	socket_type	= tcp
	protocol	= stream
	wait		= no
	user		= nobody
	group		= tty
	server		= /usr/sbin/in.ftpd
	server_args	=
	per_source	= 16
	log_type	= FILE /var/log/xinetd/ftpd
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
