# trunk/grimoire/init-scripts.d/init.d.on/bootlogd
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Default-Stop: 0 6
# Required-Start: pts ptmx
# Required-Stop: $local_fs kill
# Short-Description: runs console message capture logging daemon
### END INIT INFO

if [ "$1" == try-restart ] || ! [ -s /sbin/bootlogd ]; then exit 5; fi

. /lib/lsb/init-functions

server /sbin/bootlogd

start(){
 if   [ -d /home ] || grep -q boot=IR /proc/cmdline
 then options -l /var/log/console
 else options -l /media/root/var/log/console
 fi
}
