# trunk/grimoire/init-scripts.d/init.d.on/bootlogd
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: pts ptmx
# Default-Start: S
# Short-Description: runs console message capture logging daemon
### END INIT INFO

if [ "$1" == try-restart ] || ! [ -s /sbin/bootlogd ]; then exit; fi

. /lib/lsb/init-functions

server /sbin/bootlogd

if   [ -d /home ]
then options -l /var/log/console
else options -l /media/root/var/log/console
fi

parse "$@"
