# trunk/grimoire/init-scripts.d/init.d.on/cgroup-init
#!/bin/bash

# Copyright 2011 by Kyle Sallee,
# All rights reserved.
# For use with the Sorcerer distribution only.

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: cgroup
# Default-Start: S
# Short-Description: joins init to a cgroup
### END INIT INFO

[ $1 == start ] || exit 0

. /lib/lsb/init-functions

find_cgroup(){
 local JUNK DIR TYPE
 while read JUNK DIR TYPE JUNK; do
  if [ "$TYPE" == cgroup ] && [ -f $DIR/tasks ]; then
   CG=$DIR
   return 0
  fi
 done < /proc/mounts
 return 1
}

join(){
 find_cgroup &&
 mkdir -p $CG/daemon/init &&
 echo 1 > $CG/daemon/init/tasks
}

if   log_warning_msg "init joining cgroup"; join
then log_success_msg "init joined  cgroup"
else log_failure_msg "init failed to join cgroup"
fi
