# trunk/grimoire/init-scripts.d/init.d.on/cgroup
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Should-Start: dev
# Default-Start: S 1 2 3 4 5
# Short-Description: cgroup mounts cgroup on /dev/cgroup
# Long-Description: Read /usr/src/linux/Documentation/cgroups/cgroups.txt
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /dev/cgroup
options '-t cgroup -o blkio,clone_children,cpu,cpuset,freezer,memory cgroup'
