# trunk/grimoire/init-scripts.d/init.d.on/config
#!/bin/bash

### BEGIN INIT INFO
# Default-Start: S
# Required-Start: sys
# Short-Description: config mounts configfs on /config
### END INIT INFO

. /lib/lsb/init-functions

NAME=/config
OPTS='-n -t configfs configfs'

parse "$@"
