# trunk/grimoire/init-scripts.d/init.d.on/config.gz
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Should-Start: proc
# Default-Start: 1 2 3 4 5
# Short-Description: config.gz chmod 600 /proc/config.gz
### END INIT INFO

. /lib/lsb/init-functions

NAME=/proc/config.gz

start(){
 [ -f $NAME ] || return
 if   log_warning_msg "$NAME changing mode to 600"; chmod 600 $NAME
 then log_success_msg "$NAME mode changed"
 else log_failure_msg "$NAME mode not changed"
 fi
}

try_restart(){ :; }
parse "$@"
