# trunk/grimoire/init-scripts.d/init.d.on/console
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 1 2 3 4 5
# Short-Description: creates file /var/log/console
### END INIT INFO

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

start(){ 
 run(){
  log_warning_msg "/var/log/console creating"
  echo -n >       "/var/log/console"
  log_success_msg "/var/log/console created"
 }

 [ -f /var/log/console ] || run
}
