# trunk/grimoire/init-scripts.d/init.d.on/depmod
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: runs depmod during boot
### END INIT INFO

if [ "$1" == try-restart ] || ! [ -s /sbin/depmod ]; then exit; fi

. /lib/lsb/init-functions

log_warning_msg "depmod -A"; depmod -A
log_success_msg "depmod -A complete"
