# trunk/grimoire/init-scripts.d/init.d.on/depmod
#!/bin/bash
### BEGIN INIT INFO
# Provides: depmod
# Required-Stop: $local_fs udev
# Should-Stop: initramfs root
# Default-Stop: 0 6
# Short-Description: runs depmod as required
### END INIT INFO

# Copyright 2008 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

if [[ $1 == try-restart ]] || ! [[ -s /sbin/depmod ]]
then exit
fi

. /lib/lsb/init-functions

check(){
 VER="$( find /boot -mindepth 2 -maxdepth 2 -type f -name linux -printf "%T@	%h\n" |
         sort -gr | cut -d / -f3 | head -n 1 )"
 find /lib/modules/$VER/kernel -type f -newer \
      /lib/modules/$VER/modules.dep | grep -q .
}

if   log_warning_msg 'depmod checking modules.dep'; check
then log_warning_msg 'depmod updating'; depmod -aA $VER
     log_success_msg 'depmod updated'
fi;  log_success_msg 'depmod checked'
