# trunk/grimoire/init-scripts.d/init.d.on/empty
#!/bin/dash

### BEGIN INIT INFO
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: empty mounts tmpfs on /var/empty
### END INIT INFO

. /lib/lsb/init-functions

NAME=/var/empty
OPTS='-o defaults,size=8k,nr_inodes=4,mode=700 -t tmpfs tmpfs'

parse "$@"
