# trunk/grimoire/init-scripts.d/init.d.on/empty
#!/bin/bash

### BEGIN INIT INFO
# Provides: empty
# Required-Start: $local_fs
# Should-Start:
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: empty mounts tmpfs on /var/empty
# Description: empty mounts tmpfs on /var/empty
### END INIT INFO

. /lib/lsb/init-functions

NAME='/var/empty'
OPTS='-o defaults,size=8k,nr_inodes=4,mode=0700 -t tmpfs tmpfs'

parse "$@"
