# trunk/grimoire/init-scripts.d/init.d.on/empty
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: empty mounts tmpfs on /var/empty
### END INIT INFO

. /lib/lsb/init-functions

NAME=/var/empty
OPTS='-o defaults,size=64k,nr_inodes=16,mode=700 -t tmpfs tmpfs'

parse "$@"
