# trunk/grimoire/init-scripts.d/init.d.on/events
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: 1 2 3 4 5
# Short-Description: eliminates slugish keyboard pointer event handling during high load
### END INIT INFO

if [ "$1" = try-restart ]; then exit; fi

. /lib/lsb/init-functions

pgrep -u root events/0  | grep -q . &&
if   log_warning_msg "renice -19 \$( pgrep -u root events/0 ) executing"
                      renice -19  $( pgrep -u root events/0 ) &>/dev/null
then log_success_msg "renice -19 \$( pgrep -u root events/0 ) executed"
else log_failure_msg "renice -19 \$( pgrep -u root events/0 ) failed"
fi
