# trunk/grimoire/init-scripts.d/init.d.on/fstab
#!/bin/bash

### BEGIN INIT INFO
# Provides: $local_fs
# Required-Start: evms_activate luks proc sys
# Should-Start: cryptmount-early devices dmsetup udev vgchange
# Required-Stop: halt reboot
# Default-Start: S
# Default-Stop:
# Short-Description: fstab avoids rdinit=/bin/bash if / not in fstab
### END INIT INFO

. /lib/lsb/init-functions

start(){
 root_menu(){
  block(){ find /dev -type b -printf '%p	%Y\n' | sed '/loop/d'; }

  TITLE=$"Root Device Menu"
   HELP=$"Expected / in /etc/fstab"
 DIALOG='dialog --backtitle fstab --stdout --no-cancel'

  $DIALOG --title "$TITLE" --menu "$HELP" 0 0 0 /sbin/sulogin '' $( block )
 }

 interact(){
  ROOT=$( root_menu )
  case $ROOT in
   /sbin/sulogin) sulogin -t 300 &>/dev/console 0</dev/console ;;
   /dev*) echo "$ROOT / auto defaults 0 0" >> /etc/fstab
  esac
 }

 check(){ [ -n "$( sed 's/#.*//;s:\t: :g;\: / :p;d' /etc/fstab )" ]; }

 if   log_warning_msg 'Checking /etc/fstab for /'; check
 then log_success_msg 'Found /'
 else log_failure_msg '/ not in /etc/fstab'; interact
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
start) start ;;
try-restart) ;;
*) help ;;
esac
