# trunk/grimoire/init-scripts.d/init.d.on/fstab
#!/bin/bash

### BEGIN INIT INFO
# Should-Start: find_ir udevd
# Required-Stop: halt reboot
# Default-Start: S
# Short-Description: fstab avoids rdinit=/bin/bash if / not in fstab
# Description: fstab checks for / in /etc/fstab and rewrites it to /media/root
### END INIT INFO

. /lib/lsb/init-functions

start(){
 root_menu(){
  block(){ find /dev -type b -printf '%p	%Y\n' | sed '/loop/d'; }

  dialog --backtitle fstab --stdout --no-cancel --title "Root Device Menu" \
         --menu "Expected / in /etc/fstab" 0 0 0 /sbin/sulogin '' $( block )
 }

 interact(){
  # usb disk device may not be immediately available?
  sleep 5
  ROOT=$( root_menu )
  case $ROOT in
   /sbin/sulogin) sulogin -t 300 &>/dev/console 0</dev/console ;;
   /dev*) echo "$ROOT /media/root auto defaults 0 0" >> /etc/fstab
  esac
 }

 check(){ [ -n "$( sed 's/#.*//;s:\t: :g;\: / :p;d' /etc/fstab )" ]; }
 modify(){ sed -i 's:\t: :g;  s: / : /media/root :' /etc/fstab; }
 back(){ cp -a /etc/fstab /etc/fstab.ori; }

 [[ -d /opt ]] || grep -q boot=IR /proc/cmdline ||
 if   log_warning_msg 'Checking /etc/fstab for /'; check
 then log_success_msg 'Found /'; back; modify
 else log_failure_msg '/ not in /etc/fstab'; interact
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
start) start ;;
try-restart) ;;
*) help ;;
esac
