# trunk/grimoire/init-scripts.d/init.d.on/gconf
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: $local_fs
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: gconf updates /etc/gconf/gconf.xml.defaults based upon the content of /etc/gconf/schemas/
### END INIT INFO

# Copyright 2011 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

 DEF=/etc/gconf/gconf.xml.defaults
 SCH=/etc/gconf/schemas

check(){
 if [ -d $SCH ]; then
  ! [ -d              $DEF ] ||
  find   $SCH -cnewer $DEF | grep -q .
 else false
 fi
}

update(){
 GCONF_CONFIG_SOURCE=xml:merged:$DEF \
 gconftool-2 --makefile-install-rule $SCH/*.schemas
 for REPLY in $SCH/*.entries; do
  if   [ -f "$REPLY" ]
  then gconftool-2 --direct --config-source=xml:merged:$DEF --load $REPLY
  fi
 done
 chmod 755 $DEF
}

stop(){
 if   log_warning_msg '/etc/gconf/schemas checking'; check
 then log_warning_msg '/etc/gconf/gconf.xml.defaults updating'; update
      log_success_msg '/etc/gconf/gconf.xml.defaults updated'
 fi;  log_success_msg '/etc/gconf/schemas checked'
}
