# trunk/grimoire/init-scripts.d/init.d.on/gtk-update-icon-cache
#!/bin/dash
### BEGIN INIT INFO
# Required-Stop: kexec $local_fs root-ro udev
# Should-Stop: tmp
# Default-Stop: 0 6
# Short-Description: Updates the icon cache
### END INIT INFO

HI=/usr/share/icons/hicolor
HC=/usr/share/icons/hicolor/icon-theme.cache

up(){
 if   log_warning_msg "$HC Updating"; gtk-update-icon-cache -f -t $HI
 then log_success_msg "$HC Updated"
 else log_failure_msg "Failed to update $HC"
 fi
}

[ -f /usr/bin/gtk-update-icon-cache ] || return

. /lib/lsb/init-functions

log_warning_msg "$HC Checking"

  if ! [ -f $HC ]; then up
elif find $HI -type f -newer $HC | grep -q .; then up
else log_success_msg "$HC Current"
fi
