# trunk/grimoire/init-scripts.d/init.d.on/halt
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 0
# Required-Start: finale
# Short-Description: powers down the box
### END INIT INFO

. /lib/lsb/init-functions

deny try-restart

start(){ log_success_msg Halting; halt -d -f -i -p; }
