# trunk/grimoire/init-scripts.d/init.d.on/lilo
#!/bin/dash
### BEGIN INIT INFO
# Required-Stop: $local_fs proc udev
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: lilo regenerates /etc/lilo/conf when necessary
### END INIT INFO

# Copyright 2008-2010 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

. /lib/lsb/init-functions

NAME=/etc/lilo/conf
 CON=/etc/lilo/conf
 FAS=/etc/lilo/fast
HEAD=/etc/lilo/head
TAIL=/etc/lilo/tail

lilo_conf_update(){
 img(){
  echo "image			=	/boot/$REPLY/linux"; [ -f /etc/init.d/initramfs ] &&
  echo "	initrd		=	/boot/$REPLY/initrf"
  echo "	label		=	$REPLY"
  echo "	read-only"
  echo "	$1"
  echo
 }

 local ICX=0
 local LIM=8

 if ! [ -f $HEAD ]; then exit 1; fi

 if   [ -f /etc/init.d/initramfs ]
 then LIM=4
 fi

 cp $HEAD $CON
 cp $HEAD $FAS

 REPLY=fast img restricted >> $FAS

 find /boot -mindepth 2 -maxdepth 2 -type f -name linux -printf '%T@\t%h\n' |
 sort -gr |
 sed '\,/boot/fast,d; \,/boot/rescue,d; s:.*/::' |
 while read REPLY; do
  ICX=$(($ICX+1))
  if [ $ICX -lt $LIM ]
  then img restricted
  else
   rm -fr /boot/$REPLY \
   /lib/modules/$REPLY \
   /lib/modules/$REPLY.old
  fi
 done | tee -a $FAS >> $CON

 if   [ -f /boot/rescue/initrf ]
 then REPLY=rescue img mandatory | tee -a $FAS >> $CON
 fi

 if   [ -f $TAIL ]
 then cat  $TAIL | tee -a $FAS >> $CON
 fi

 if [ -x /etc/init.d/initramfs ]
 then    /etc/init.d/initramfs
 else /sbin/lilo -P ignore -C $CON
 fi
}

both(){
 find /boot -mindepth 2 -maxdepth 2 -type f -name linux -printf '%P\n'
 sed 's:#.*::;s: ::g;s:\t::g;s:image=/boot/::p;d' $CON
}

current(){ 
 [ $CON -nt $HEAD ] &&
 [ $CON -nt $TAIL ] &&
 ! both | sed '\,fast/linux,d' | sort | uniq -u | grep -q .
}

begin(){
 if   current
 then return
 elif log_warning_msg "$NAME updating"; lilo_conf_update
 then log_success_msg "$NAME updated"
 else log_failure_msg "$NAME failed"
  if   [ -n $1 ]
  then sulogin -t 300 1>/dev/console 2>&1 0</dev/console
  fi
 fi
}

main(){
 if ! [ -d /boot ] ||
      [ -z "$( find /boot -mindepth 2 -maxdepth 2 -type f -name linux )" ]
 then return
 fi
 log_warning_msg "$NAME checking"; begin $1
 log_success_msg "$NAME checked"
}

export LC_ALL=C
case $1 in
 stop) main s ;;
    *) main ;;
esac
