# trunk/grimoire/init-scripts.d/init.d.on/lock
#!/bin/bash

### BEGIN INIT INFO
# Provides: lock
# Required-Start: proc
# Should-Start:
# Required-Stop: 
# Default-Start: S 1 2 3 4 5
# Default-Stop:
# Short-Description: lock mounts tmpfs on /var/lock
# Description: lock mounts tmpfs on /var/lock
### END INIT INFO

. /lib/lsb/init-functions

NAME='/var/lock'
OPTS='-o defaults,size=1m,nr_inodes=1k,mode=0755 -t tmpfs tmpfs'

parse "$@"
