# trunk/grimoire/init-scripts.d/init.d.on/loglevel
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: S
# Short-Description: loglevel sets the level at which the kernel will write to the console.
# Description: The kernel routine printk() will only print a message on the console,
#	if the message loglevel value is less than
#	the value of the variable console_loglevel.
#	For a description of the values for console_loglevel execute: man syslog
### END INIT INFO

. /lib/lsb/init-functions

start(){
 if   log_warning_msg 'Kernel logging to console level 3 setting'; echo 3 > /proc/sys/kernel/printk
 then log_success_msg 'Loglevel 3 set'
 else log_failure_msg 'Loglevel 3 not set'
 fi
}
       stop(){ :; }
    restart(){ :; }
try_restart(){ :; }
parse "$@"
