# trunk/grimoire/init-scripts.d/init.d.on/modprobe
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc sys modprobe-coldplug modprobe-r udev
# Default-Start: S
# Short-Description: modprobe loads modules during boot
### END INIT INFO

if [ "$1" == try-restart ]; then exit 5; fi

. /lib/lsb/init-functions

if [ -x /sbin/modprobe ] &&
   [ -f /etc/modules   ]
then
 sed 's:#.*::;/^ *$/d' /etc/modules |
 while read REPLY; do
  if   log_warning_msg "Module $REPLY loading" && modprobe $REPLY
  then log_success_msg "Module $REPLY loaded"
  else log_failure_msg "Module $REPLY not loaded"
  fi
 done
fi
