# trunk/grimoire/init-scripts.d/init.d.on/modprobe-r
#!/bin/dash

### BEGIN INIT INFO
# Required-Start: proc sys
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: modprobe-r unloads modules during boot
# Description: modprobe-r unloads modules specified in /etc/modules.rm during boot
### END INIT INFO

if [ "$1" = try-restart ] || ! [ -f /etc/modules.rm ] || ! [ -x /sbin/modprobe ]; then exit; fi

. /lib/lsb/init-functions

rmp(){ modprobe -r $( sed 's:#.*::;/^ *$/d' /etc/modules.rm ); }

log_warning_msg 'Modules unloading'; rmp
log_success_msg 'Modules unloaded'
