# trunk/grimoire/init-scripts.d/init.d.on/modprobe-r
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Required-Start: modprobe-coldplug proc sys
# Should-Start: udev
# Default-Start: S
# Short-Description: modprobe-r unloads modules specified in /etc/modules.rm during boot
### END INIT INFO

if [ "$1" == try-restart ] || ! [ -f /etc/modules.rm ] || ! [ -x /sbin/modprobe ]; then exit; fi

. /lib/lsb/init-functions

rmp(){
 [ -f /etc/modules.rm ] &&
 sed 's:#.*::;/^ *$/d' /etc/modules.rm |
 xargs -r --max-lines=64 modprobe -r
}

log_warning_msg 'Modules unloading'; rmp
log_success_msg 'Modules unloaded'
