# trunk/grimoire/init-scripts.d/init.d.on/mount
#!/bin/dash

### BEGIN INIT INFO
# Provides: $local_fs
# Required-Start: mtab root-rw move_dev
# Required-Stop: cryptmount cryptsetup halt mdadm reboot reboot-kexec
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mount mounts and un-mounts local file systems
### END INIT INFO

. /lib/lsb/init-functions

start(){
 run(){
  if [ -f     /media/root/sbin/init ]
  then chroot /media/root mount -a
  else                    mount -a
  fi
 }

 if   log_warning_msg 'Filesystems from /etc/fstab mounting'; run
 then log_success_msg 'Filesystems from /etc/fstab mounted'
 else log_success_msg 'Filesystems from /etc/fstab not mounted'
 fi
}

stop(){
 run(){
  sync
  [ -f /proc/mounts ] || mount -nrt proc proc /proc

  cut -d ' ' -f2 /proc/mounts |
  LC_ALL=C sort -ur | xargs -r --max-lines=1024 umount -lr
 }

 log_warning_msg 'Filesystems in /proc/mounts unmounting or remounting read-only'; run
 log_success_msg 'Filesystems in /proc/mounts unmounted  or remounted  read-only'
}

try_restart(){ :; }
parse "$@"
