# trunk/grimoire/init-scripts.d/init.d.on/mount
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Provides: $local_fs
# Required-Start: mtab root-rw bind_dev
# Required-Stop: cryptmount cryptsetup finale halt mdadm reboot reboot-kexec root-ro vgchange
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mount mounts and un-mounts local file systems
### END INIT INFO

. /lib/lsb/init-functions

start(){
 run(){
  if [ -f     /media/root/sbin/init ]
  then chroot /media/root mount -a
  else                    mount -a
  fi
 }

 if   log_warning_msg 'Filesystems from /etc/fstab mounting'; run
 then log_success_msg 'Filesystems from /etc/fstab mounted'
 else log_success_msg 'Filesystems from /etc/fstab not mounted'
 fi
}

stop(){
 run(){
  sync
  [ -f /proc/mounts ] || mount -nrt proc proc /proc

   cut -d ' ' -f2 /proc/mounts |
   LC_ALL=C sort -r |
   if   [ -d /media/root/proc ]
   then grep /media/root
   else grep -vx / ; mount -o remount,ro /
   fi |
   xargs -r --max-lines=1024 umount -lr
 }

 log_warning_msg 'Filesystems unmounting'; run
 log_success_msg 'Filesystems unmounted'
}

try_restart(){ :; }
parse "$@"
