# trunk/grimoire/init-scripts.d/init.d.on/move
#!/bin/bash

### BEGIN INIT INFO
# Provides: move
# Required-Start: $local_fs
# Default-Start: S
# Short-Description: move moves mounted virtual filesystems into /media/root
### END INIT INFO

. /lib/lsb/init-functions

start(){
 run(){
  sed '\,rootfs,d
       \,/media/root,d
       \,^/dev/,d
       \,^/sys/,d' /proc/mounts |
  while read DEV MNT JNK; do
   [[ -d        $MNT ]] || continue
   [[ -d             /media/root$MNT ]] ||
   mkdir -p          /media/root$MNT &&
   mount --move $MNT /media/root$MNT
  done
 }

 [[ -f /media/root/sbin/init ]] || return
 log_warning_msg 'Moving mounted filesystems to /media/root'; run
 log_success_msg 'Moving mounted filesystems to /media/root'
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case "$@" in
 start) start ;;
 try-restart) : ;;
 *) help ;;
esac
