# trunk/grimoire/init-scripts.d/init.d.on/proc
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: S
# Short-Description: proc mounts proc on /proc
### END INIT INFO

. /lib/lsb/init-functions

only start
name /proc
options '-n -t proc proc'
