# trunk/grimoire/init-scripts.d/init.d.on/proc
#!/bin/bash

### BEGIN INIT INFO
# Provides: proc
# Required-Start: 
# Should-Start: 
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: proc mounts proc on /proc
# Description: proc mounts proc on /proc
### END INIT INFO

. /lib/lsb/init-functions

SERV='mount'
NAME='/proc'
OPTS='-n -t proc proc'

parse "$@"
