# trunk/grimoire/init-scripts.d/init.d.on/ptmx
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: pts
# Default-Start: S 1 2 3 4 5
# Short-Description: creates symbolic link from pts/ptmx to /dev/ptmx
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 if   [ -c       /dev/pts/ptmx ] &&
    ! [ -h       /dev/ptmx     ]; then
  rm -f          /dev/ptmx
  ln -s pts/ptmx /dev/ptmx
 fi
}
