# trunk/grimoire/init-scripts.d/init.d.on/pts
#!/bin/bash

### BEGIN INIT INFO
# Provides: pts
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: pts mounts devpts on /dev/pts
### END INIT INFO

. /lib/lsb/init-functions

NAME=/dev/pts
OPTS='-t devpts devpts'

parse "$@"
