# trunk/grimoire/init-scripts.d/init.d.on/pts
#!/bin/bash

### BEGIN INIT INFO
# Provides: pts
# Required-Start: $local_fs
# Should-Start: udev
# Required-Stop: 
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: devpts mounts devpts on /dev/pts
# Description: devpts mounts devpts on /dev/pts
### END INIT INFO

. /lib/lsb/init-functions


start()  {
  if    log_warning_msg  '/dev/pts mounting'; mount -t devpts devpts /dev/pts
  then  log_success_msg  '/dev/pts mounted'
  else  log_failure_msg  '/dev/pts not mounted'
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
