# trunk/grimoire/init-scripts.d/init.d.on/reboot
#!/bin/dash

### BEGIN INIT INFO
# Default-Stop: 6
# Short-Description: reboot kills lingering processes and reboots the box
### END INIT INFO

. /lib/lsb/init-functions

case $1 in
stop) log_success_msg Rebooting; reboot -d -f -i  ;;
esac
