# trunk/grimoire/init-scripts.d/init.d.on/reboot
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 6
# Required-Start: finale
# Should-Start: reboot-kexec
# Short-Description: reboots the box
### END INIT INFO

. /lib/lsb/init-functions

deny try-restart

start(){ log_success_msg Rebooting; reboot -d -f -i ; }
