# trunk/grimoire/init-scripts.d/init.d.on/root-rw
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: fsck
# Required-Stop: $local_fs
# Default-Start: S
# Short-Description: At boot mounts / or /media/root rw
### END INIT INFO

. /lib/lsb/init-functions

start(){
 run(){
  if [ -d /opt ]
  then mount -n -o remount,rw /
  else mount -n -o         rw /media/root
  fi
 }

 if   log_warning_msg 'root mounting read write'; run
 then log_success_msg 'root mounted  read write'
 else log_failure_msg 'root not mounted'
 fi
}

try_restart(){ :; }
stop(){ :; }
parse "$@"
