# trunk/grimoire/init-scripts.d/init.d.on/rtc-cmos
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: udev
# Default-Start: S
# Short-Description: rtc-cmos loads the module
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
