# trunk/grimoire/init-scripts.d/init.d.on/run
#!/bin/bash

### BEGIN INIT INFO
# Provides: run
# Required-Start: proc $local_fs
# Should-Start:
# Required-Stop: 
# Default-Start: S 
# Default-Stop:
# Short-Description: run mounts tmpfs on /var/run
# Description: run mounts tmpfs on /var/run
### END INIT INFO

. /lib/lsb/init-functions

NAME='/var/run'
OPTS='-o defaults,size=1m,nr_inodes=1k,mode=0755 -t tmpfs tmpfs'

parse "$@"
