# trunk/grimoire/init-scripts.d/init.d.on/run
#!/bin/bash

### BEGIN INIT INFO
# Provides: run
# Required-Start: $local_fs
# Should-Start:
# Required-Stop: 
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: run mounts tmpfs on /var/run
# Description: run mounts tmpfs on /var/run
### END INIT INFO

. /lib/lsb/init-functions


OPTIONS='defaults,size=1m,nr_inodes=1k,mode=0750'

start()  {
  run()  {
    if    !  [   -d  /var/run  ]
    then  mkdir  -p  /var/run
    fi  &&
    mount  -o $OPTIONS  -t tmpfs  tmpfs  /var/run
  }

  if    ! sed   's/#.*//'   /etc/fstab  |
          grep  -q          /var/run;  then
    if    log_warning_msg  '/var/run mounting';  run
    then  log_success_msg  '/var/run mounted'
    else  log_failure_msg  '/var/run not mounted'
    fi
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
