# trunk/grimoire/init-scripts.d/init.d.on/run
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: S 1 2 3 4 5
# Short-Description: run mounts tmpfs on /run
### END INIT INFO

. /lib/lsb/init-functions

NAME=/run
OPTS='-o defaults,size=1m,nr_inodes=2k,mode=0755 -t tmpfs tmpfs'

parse "$@"
