# trunk/grimoire/init-scripts.d/init.d.on/scsi_wait_scan
#!/bin/dash

### BEGIN INIT INFO
# Should-Start: udevsettle
# Default-Start: S
# Short-Description: scsi_wait_scan loads the module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
