# trunk/grimoire/init-scripts.d/init.d.on/services
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: $local_fs
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: services updates /etc/services based upon the content of /etc/services.d/
### END INIT INFO

# Copyright 2011 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

if [ "$1" == try-restart ]; then exit; fi

. /lib/lsb/init-functions

check(){
 ! [ -f /etc/services ] ||
 find   /etc/services.d -type f -cnewer /etc/services | grep -q .
}

update(){
 find /etc/services.d/ -maxdepth 1 -type f |
 xargs -r --max-lines=256 \
 sed 's: *#.*::;s:  *:\t:g' > /etc/services
}

if   log_warning_msg '/etc/services.d checking'; check
then log_warning_msg '/etc/services   updating'; update
     log_success_msg '/etc/services   updated'
fi;  log_success_msg '/etc/services.d checked'
