# trunk/grimoire/init-scripts.d/init.d.on/shells
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: $local_fs
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: shells updates /etc/shells based upon the content of /etc/shells.d/
### END INIT INFO

# Copyright 2011 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

check(){
 ! [ -f /etc/shells ] ||
 find /etc/shells.d -cnewer /etc/shells | grep -q .
}

update(){
 find /etc/shells.d/ -type l -printf "%l\n" > /etc/shells
}

stop(){
 if   log_warning_msg '/etc/shells.d checking'; check
 then log_warning_msg '/etc/shells   updating'; update
      log_success_msg '/etc/shells   updated'
 fi;  log_success_msg '/etc/shells.d checked'
}
