# trunk/grimoire/init-scripts.d/init.d.on/shm
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: shm mounts tmpfs on /dev/shm
### END INIT INFO

. /lib/lsb/init-functions

NAME=/dev/shm
OPTS='-t tmpfs tmpfs'

parse "$@"
