# trunk/grimoire/init-scripts.d/init.d.on/sorcery
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: sorcery mounts tmpfs on /root/.sorcery
### END INIT INFO

. /lib/lsb/init-functions

NAME=/root/.sorcery
OPTS='-o defaults,size=512m,nr_inodes=512k,mode=700 -t tmpfs tmpfs'

parse "$@"
