# trunk/grimoire/init-scripts.d/init.d.on/sys
#!/bin/bash

### BEGIN INIT INFO
# Provides: sys
# Required-Start: 
# Should-Start: 
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: sys mounts sysfs on /sys
# Description: sys mounts sysfs on /sys
### END INIT INFO

. /lib/lsb/init-functions

NAME='/sys'
OPTS='-n -t sysfs sysfs'

parse "$@"
