# trunk/grimoire/init-scripts.d/init.d.on/sys
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Default-Start: S
# Required-Start: proc
# Short-Description: sys mounts sysfs on /sys
### END INIT INFO

. /lib/lsb/init-functions

NAME=/sys
OPTS='-n -t sysfs sysfs'

parse "$@"
