# trunk/grimoire/init-scripts.d/init.d.on/sys
#!/bin/bash

### BEGIN INIT INFO
# Provides: sys
# Required-Start: 
# Should-Start: 
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: sys mounts sysfs on /sys
# Description: sys mounts sysfs on /sys
### END INIT INFO

. /lib/lsb/init-functions


start()  {
  run()  {
    if    !  [   -d  /sys  ]
    then  mkdir  -p  /sys
    fi  &&
    mount  -o $OPTIONS  -n  -t sysfs  sysfs  /sys
  }

  if    [  -e  /sys/bus  ]
  then  log_success_msg  '/sys already mounted'
  elif  log_warning_msg  '/sys mounting';  run
  then  log_success_msg  '/sys mounted'
  else  log_failure_msg  '/sys not mounted'
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
