# trunk/grimoire/init-scripts.d/init.d.on/tmp
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Required-Start: $local_fs
# Required-Stop: proc swap
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts a 3G tmpfs on /tmp
### END INIT INFO

. /lib/lsb/init-functions

NAME=/tmp
OPTS='-o defaults,size=3g,nr_inodes=1m,mode=1777 -t tmpfs tmpfs'
parse "$@"
