# trunk/grimoire/init-scripts.d/init.d.on/update-mime-database
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Stop: kexec $local_fs root-ro udev
# Should-Stop: tmp
# Default-Stop: 0 6
# Short-Description: Updates mime database in mime directories
### END INIT INFO


check(){
 local DIR
 while read DIR; do
  log_warning_msg "$DIR checking"
    if ! [ -d $DIR ]; then :
  elif ! [ -d $DIR/packages ]; then :
  elif ! [ -f $DIR/mime.cache ]; then echo "$DIR"
  elif find $DIR/packages -maxdepth 0 -newer $DIR/mime.cache | grep -q  .; then echo "$DIR"
  fi
  log_success_msg "$DIR checked"
 done
}

run(){
 while read DIR; do
  if   log_warning_msg "update-mime-database $DIR executing"; update-mime-database $DIR
  then log_success_msg "update-mime-database $DIR executed"
  else log_failure_msg "update-mime-database $DIR failed"
  fi
 done
}

[ -f /usr/bin/update-mime-database ] || exit 0

. /lib/lsb/init-functions

stop(){
 echo "/usr/share/mime
/opt/kde4/share/mime
/opt/kde3/share/mime" |
 check | run
}
