# trunk/grimoire/init-scripts.d/init.d.on/usr_src_sorcery_tmpfs
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs
# Required-Stop: proc swap
# Should-Start: usr_src_sorcery_bind
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mounts 3G tmpfs on /usr/src/sorcery
### END INIT INFO

RAM=$( sed -r "s/MemTotal: *(.*) kB/\1/p;d" /proc/meminfo )
if (( 1048576 > RAM )); then exit 0 ; fi

. /lib/lsb/init-functions


NAME=/usr/src/sorcery
OPTS='-o defaults,size=4g,nr_inodes=1m,mode=700 -t tmpfs tmpfs'
parse "$@"
