# trunk/grimoire/init-scripts.d/init.d.on/var_lock_run
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: lock
# Default-Start: S 1 2 3 4 5
# Short-Description: creates symbolic links /var/{lock,run}
### END INIT INFO

if [ "$1" == try-restart ]; then exit 5; fi

. /lib/lsb/init-functions

start(){
 if ! [ -d /run/lock ]; then mkdir -p /run/lock; fi
 if ! [ -d /var      ]; then mkdir    /var; fi
 if ! [ -h /var/lock ]; then rm -fr /var/lock; ln -s /run/lock /var/lock; fi
 if ! [ -h /var/run  ]; then rm -fr /var/run;  ln -s /run      /var/run; fi
}
