# trunk/grimoire/kexec-tools.d/init.d.on/kexec
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Default-Stop: 6
# Required-Stop: mount
# Short-Description: kexec prepares box for reboot faster.
### END INIT INFO

. /lib/lsb/init-functions

run(){
 if   [ -f /boot/fast/initrf ]
 then
  I='--initrd=/boot/fast/initrf'
  L='-l /boot/fast/linux'
 else
  V="$( find /boot -maxdepth 2 -mindepth 2 -name linux -printf '%T@\t%h\n' |
        sort -gr | head -n 1 | sed 's:.*/::' )"
  [ -f /boot/$V/initrf ] &&
  I="--initrd=/boot/$V/initrf"
  L="-l /boot/$V/linux"
 fi

 log_warning_msg "kexec $L $I"; kexec $L $I
}

case $1 in
stop) run ;;
esac
