# trunk/grimoire/kexec-tools.d/init.d.on/reboot-kexec
#!/bin/dash

### BEGIN INIT INFO
# Default-Stop: 6
# Required-Stop: reboot
# Short-Description: reboot-kexec soft boots with previously loaded linux and initrf
### END INIT INFO

. /lib/lsb/init-functions

run(){
 log_success_msg Rebooting; kexec -e
 kexec -e
}

case $1 in
stop) run ;;
esac
