# trunk/grimoire/kexec-tools.d/init.d.on/reboot-kexec
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Start: 6
# Required-Start: finale
# Short-Description: reboot-kexec soft boots with previously loaded linux and initrf
### END INIT INFO

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

start(){
 log_success_msg Rebooting; kexec -e
 kexec -e
}
