# trunk/grimoire/networks.d/init.d.on/hostname
#!/bin/bash

### BEGIN INIT INFO
# Provides: hostname
# Required-Start: proc
# Required-Stop: 
# Default-Start: S
# Default-Stop: 
# Short-Description: hostname sets the box's name
# Description: hostname sets the box's name
### END INIT INFO

. /lib/lsb/init-functions

NAME=hostname

start(){
 log_warning_msg "$NAME setting"
 [ -f   /etc/hostname ] &&
 read < /etc/hostname
 echo ${REPLY:-noname} > /proc/sys/kernel/hostname
 log_success_msg "$NAME set"
}

parse "$@"
