# trunk/grimoire/networks.d/init.d.on/hostname
#!/bin/bash

### BEGIN INIT INFO
# Provides: hostname
# Required-Start: proc
# Required-Stop: 
# Default-Start: 1 2 3 4 5
# Default-Stop: 
# Short-Description: hostname sets the box's name
# Description: hostname sets the box's name
### END INIT INFO

. /lib/lsb/init-functions

start() {
 log_warning_msg "hostname setting"
 read        <             /etc/hostname
 echo $REPLY > /proc/sys/kernel/hostname
 log_success_msg "hostname set"
}

case $1 in
 try-restart) :     ;;
       start) start ;;
esac
