# trunk/grimoire/networks.d/init.d.on/networks
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: $network
# Required-Start: hostname loopback run var_lock_run
# Required-Stop: $local_fs
# Should-Stop: kill
# Should-Start: modprobe udev udevsettle
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: networks raises and lowers network connections
### END INIT INFO

[ "$1" == try-restart ] && exit 2

. /lib/lsb/init-functions

# Please do not modify this file.
# Modify files located in    /etc/networks.d/
# Templates are in           /etc/networks.d/templates/
# If unfamiliar please read  /etc/networks.d/templates/README

# If the networking should not start during boot then execute:
# cast -r networks

NETWORKS=/etc/networks.d

NAME=networks
LOCK=/var/lock/networks

start(){
 log_warning_msg "$NAME checking"; [ -f $LOCK ] &&
 log_success_msg "$NAME already started" && return
 log_warning_msg "$NAME starting"

 for     file in $NETWORKS/*; do
  [ -f "$file" ]    &&
  bash "$file" start &&
  break
 done &&
 echo -n > $LOCK &&
 log_success_msg "$NAME started" ||
 log_failure_msg "$NAME not started"
}

stop(){
 log_warning_msg "$NAME stopping"
 rm -f $LOCK
 for    file in $NETWORKS/*; do
 [ -f "$file" ] &&
 bash "$file" stop
 done
 log_success_msg "$NAME has stopped"
}

restart(){ stop; sleep 3; start; }
