# trunk/grimoire/avahi.d/init.d/avahi-daemon
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 510
# Provides: avahi
# Required-Start: dbus $network run
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Required-Stop: kill
# Should-Stop: kill fuse.gunzip
# Short-Description: avahi-daemon is a ZeroConf daemon for mDNS and service registration
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart reload status
server /usr/sbin/avahi-daemon
name avahi-daemon
options -D

[     -d /var/run/avahi-daemon ] ||
mkdir -p /var/run/avahi-daemon

stop(){
 if   log_warning_msg "stop""; $SERV -k
 then log_success_msg "stop"
 else log_failure_msg "stop"; return 1
 fi
}

reload(){
 if   log_warning_msg "reload"; $SERV -r
 then log_success_msg "reload"
 else log_failure_msg "reload"; return 1
 fi
}
