# trunk/grimoire/acpid.d/init.d.on/acpid
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: acpid starts the ACPI Event Daemon
# Description: acpid is the ACPI Event Daemon for power management and
#     delivers ACPI events from the kernel to user space.
### END INIT INFO

. /lib/lsb/init-functions

only start
server /usr/sbin/acpid
