# trunk/grimoire/courier.d/init.d/courier
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: $mail
# Required-Start: $network
# Should-Start: $remote_fs
# Required-Stop: $network
# Should-Stop:  $remote_fs
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: courier is a MTA
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status

run(){ if   [ -x $1 ]; then $1 $2; fi; }

block(){
 run /usr/sbin/courierfilter $1
 run /usr/libexec/authlib/authdaemond $1
 run /usr/sbin/courierldapaliasd $1
 run /usr/libexec/courier/pcpd $1
 run /usr/sbin/courier $1
 run /usr/sbin/esmtpd $1
 run /usr/sbin/esmtpd-msa $1
}  

start(){ echo "Courier MTA starting"; block start; }
 stop(){ echo "Courier MTA stopping"; block stop; }
