# trunk/grimoire/irrd.d/init.d/irrd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network
# Should-Start: $remote_fs
# Required-Stop: $network
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: Internet Routing Registry Daemon
# Description: irrd is a database for RIPE-181 and RPSL routing.
### END INIT INFO

. /lib/lsb/init-functions

only start stop status
server /sbin/irrd
config /etc/irrd.conf
options '-g irrd -l irrd'
