# trunk/grimoire/sorcery-bxz.d/init.d.on/bxz
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: bxz xz compresses tarballs in the bacground
# Description: bxz provides background compression of tarballs.
# bxz enables more disk space to be conserved by
# recompressing tarballs more efficiently in the background
# while the box has a low loadavg
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/bxz.pid
SERV=/usr/sbin/bxz

parse "$@"
