# trunk/grimoire/procps.d/init.d.on/sysctl
#!/bin/bash

### BEGIN INIT INFO
# Provides: sysctl
# Required-Start:
# Should-Start: udev modprobe
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: sysctl sets kernel parameters
# Description:       sysctl sets kernel parameters
#	configuration file is /etc/sysctl.conf 
#	when /etc/sysctl.$RUNLEVEL.conf is not present
### END INIT INFO

. /lib/lsb/init-functions

NAME=sysctl
SERV=/sbin/sysctl
CONF=/etc/sysctl.$RUNLEVEL.conf; [ -f $CONF ] ||
CONF=/etc/sysctl.conf
OPTS="-p $CONF"

parse "$@"
