# trunk/grimoire/sysvinit.d/inittab
id:3:initdefault:

# Runlevels
# 0: halt
# 1: sequential init
# 2: parallel   init + multi-user + networking
# 3: sequential init + multi-user + networking
# 4: parallel   init + multi-user + networking + GUI login
# 5: sequential init + multi-user + networking + GUI login
# 6: reboot

# Runlevels 0 and 6 are absolute.
# Every root file system has it own /etc/inittab
# with it's own initdefault
# Runlevels 3 and 5 are tested for bugs.
# Runlevels 2 and 4 might be bugged.
