# trunk/grimoire/alsa-utils.d/init.d/snd-seq-oss
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: snd-mixer-oss
# Default-Start: S
# Short-Description: snd-seq-oss loads the module
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
