# trunk/grimoire/alsa-utils.d/init.d/snd-pcm-oss
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: modprobe udev udevtrigger
# Default-Start: S
# Short-Description: snd-pcm-oss loads the module
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
