# trunk/grimoire/laptop-mode-tools.d/init.d/laptop-mode
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs $syslog run
# Required-Stop: $local_fs
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: laptop-mode extends battery life.
### END INIT INFO

# config:  /etc/laptop-mode/laptop-mode.conf

rootfs || exit 0

. /lib/lsb/init-functions

# Note: Location of PIDF is hardcoded in /usr/bin/laptop_mode!
PDIR=/var/run/laptop-mode-tools

only start stop restart try-restart status
#pid_file $PDIR/enabled
server /sbin/laptop_mode
name laptop-mode

start(){
 run(){ /bin/mkdir -p "$PDIR" && /bin/touch "$PIDF" && "$SERV" auto init; }
 if   log_warning_msg "laptop_mode starting"; run
 then log_success_msg "laptop_mode started"
 else log_failure_msg "laptop_mode start failure"
 fi
}

stop(){
 die(){ /bin/rm -f "$PIDF"; "$SERV" stop init; }

 if   log_warning_msg "laptop_mode stopping"; die
 then log_success_msg "laptop_mode stopped"
 else log_success_msg "laptop_mode stop failure"
 fi
}

restart(){ "$SERV" stop; /bin/sleep 1; "$SERV" start; }
status(){  "$SERV" status; }
