# trunk/grimoire/laptop-mode-tools.d/init.d/laptop-mode
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs $syslog run
# Required-Stop: $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: laptop-mode extends battery life.
### END INIT INFO

# config:  /etc/laptop-mode/laptop-mode.conf

. /lib/lsb/init-functions

# Note: Location of PIDF is hardcoded in /usr/bin/laptop_mode!
PDIR=/var/run/laptop-mode-tools

only start stop restart try-restart status
pid_file $PDIR/enabled
server /usr/sbin/laptop_mode
name laptop-mode

start(){ mkdir -p $PDIR; touch  $PIDF; $SERV auto init; }
stop(){ rm -f $PIDF; $SERV stop init; }
restart(){ $SERV stop; $SERV start; }
status(){ $SERV status; }    
