# trunk/grimoire/sorcerer-installer.d/init.d.on/alt_boot
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: archives
# Should-Start: cryptsetup mdadm udevsettle vgchange
# Default-Start: S
# Short-Description: boot presents additional boot options when booting from I/R
### END INIT INFO

# Copyright 2007 through 2011 by Kyle Sallee,
# All rights reserved.
# For use with the Sorcerer distribution only.

. /lib/lsb/init-functions

only start

start(){
 block(){ find /dev -type b -printf '%p	%Y\n' | sed '/loop/d'; }

 root_menu(){
  dialog --backtitle alt_boot --stdout --no-cancel --title "Root Device Menu" \
         --menu "Please select device for /" 0 0 0 $( block )
 }

 root(){ echo "$( root_menu ) /media/root auto defaults 0 0" >> /etc/fstab; }

 main_menu(){
  dialog --backtitle alt_boot --item-help --stdout --no-cancel --timeout 60 \
         --title "Boot Options Menu" \
         --menu  "Please select an option within 60 seconds" 0 0 0 \
         "IR"    "Continue booting Install Rescue and mount the Install Rescue initramfs as /" \
         "Select for installation or manual rescue" \
         "Root"   "Continue booting Install Rescue however mount an alternate /" \
         "Select to boot while avoid LiLo errors and initramfs errors"
 }

 main(){
  case $( main_menu ) in
   Root) root ;;
   *)  return ;;
  esac
 }


 grep -q boot=IR /proc/cmdline || return

 if   log_warning_msg 'boot options menu'; main
 then log_success_msg 'boot options menu complete'
 else log_failure_msg 'boot options menu incomplete'
 fi
}
