# trunk/grimoire/sorcerer-installer.d/init.d.on/archives
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: find_ir
# Should-Start: udev loadkeys
# Default-Start: S
# Short-Description: archives makes the initramfs more usable.
### END INIT INFO

. /lib/lsb/init-functions

deny try-restart

start(){
 load(){
  cd /
  cat $FROM/{coreutils,dbus,dbus-glib,expat,file,glib,inetutils,kexec-tools,metalog,openssh,openssl,parted,pciutils,polkit,procps,psmisc,usbutils,util-linux}/*/*.tar.xz |
  xz -cd | tar -ixpf -
  local -i AVAIL
  AVAIL=$( free -mt | grep Total: | tr -s ' ' | cut -d ' ' -f4 )
  if (( AVAIL > 300 )); then
   cat $FROM/{btrfs-tools,e2fsprogs,gcc,groff,jfsutils,man,man-pages,reiserfsprogs,util-linux,xfsprogs}/*/*.tar.xz |
   xz -cd | tar -ixpf -
  fi
  rm -fr /usr/include /usr/lib{,64}/*.a /usr{,/share}/info \
         /usr{,/share}/doc /usr/share/gtk-doc
  ldconfig
 }

 local FROM=/media/IR/var/cache/archive/$HOSTTYPE
 [ -d $FROM ] &&
 if   log_warning_msg 'archives loading'; load
 then log_success_msg 'archives loaded'
 else log_failure_msg 'archives not loaded'
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }
