# trunk/grimoire/sorcerer-installer.d/init.d.on/archives
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: find_ir
# Should-Start: udev loadkeys
# Default-Start: S
# Short-Description: archives makes the initramfs more usable.
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 load(){
  cd /
  cat $FROM/{bash,btrfs-tools,coreutils,cryptsetup,dbus,dbus-glib,e2fsprogs,expat,extlinux,fhs,file,gcc,glib,groff,inetutils,jfsutils,kexec-tools,man,man-pages,mdadm,metalog,nano,openssh,openssl,parted,pciutils,polkit,procps,psmisc,reiserfsprogs,sysvinit,usbutils,util-linux,xfsprogs}/*/*.tar.xz |
  xz -cd | tar -ixpf -
 
  rm -fr /usr/include /usr/lib{,64}/*.a /usr/share/gtk-doc
  ldconfig
 }

 local FROM=/media/IR/var/cache/archive/$HOSTTYPE
 [ -d $FROM ] &&
 if   log_warning_msg 'archives loading'; load
 then log_success_msg 'archives loaded'
 else log_failure_msg 'archives not loaded'
 fi
}
