# trunk/grimoire/sorcerer-installer.d/init.d.on/welcome
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: find_ir $local_fs move udev
# Default-Start: S
# Short-Description: welcome instructs the SA how to begin installation
### END INIT INFO

. /lib/lsb/init-functions

start(){
 on_ir(){ grep -q boot=IR /proc/cmdline; }

 welcome(){
  echo
  echo $"Welcome to Sorcerer."
  echo $"Please login as root."
  echo $"Press enter for password."
  echo $"/sbin/sinstall"
  echo $"to begin a menu driven installation."
  echo
  echo $"Or if one is a new SA, system administrator,"
  echo $"then please first complete the dialog tutorial."
  echo $"/sbin/dt"
  echo $"to begin the dialog tutorial."
 }

 if   on_ir &&
      log_warning_msg "welcome starting"; welcome > /dev/console
 then log_success_msg "welcome completed"
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
       start) start ;;
 try-restart) :     ;;
           *) help  ;;
esac
