# trunk/grimoire/xen.d/init.d/xen
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $syslog $remote_fs proc
# Required-Stop:  $syslog $remote_fs
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mount xenfs on /proc/xen
### END INIT INFO

[ -d /proc/xen ] || exit 0

. /lib/lsb/init-functions
only start stop restart status
name /proc/xen
options '-t xenfs xenfs /proc/xen'
