# trunk/grimoire/xen.d/init.d/xen-watchdog
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $syslog $remote_fs
# Should-Start: xend
# Required-Stop: $syslog $remote_fs
# Should-Stop: xend
# Default-Start: 1 2 3 5
# Default-Stop: 0 6
# Short-Description: starts and stops xen-watchdog
### END INIT INFO

# Probably should run on only one root file system.

. /lib/lsb/init-functions
server /usr/sbin/xenwatchdogd
options "30 15"
only start stop restart try-restart status
