# trunk/grimoire/immune-sentinel-resume.d/init.d.on/sentinel-resume
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: resumes processing on delayed and paused spells
### END INIT INFO

. /lib/lsb/init-functions

pid_file /var/run/sentinel-resume.pid
server   /sbin/sentinel-resume

parse "$@"
