# trunk/grimoire/bind.d/init.d.on/capability
#!/bin/bash

### BEGIN INIT INFO
# Provides: capability
# Should-Start: udev
# Required-Start:
# Required-Stop:
# Default-Start: 3 4 5
# Default-Stop:
# Short-Description: capability loads the module
# Description: capability loads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME=capability

parse "$@"
