# trunk/grimoire/snort.d/init.d.on/snort
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $network
# Required-Stop:  $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: snort provides an intrusion detection and response system.
### END INIT INFO

# incomplete init file.
# Must get new snort up and running to see what pid file
# and where it puts it, and if automatic, and complete
# the functions for restarting, reloading, and stopping snort.

. /lib/lsb/init-functions

config /etc/snort/snort.conf
server /usr/bin/snort
name snort

start(){
 if   log_warning_msg "$NAME starting"; start_daemon -p "$PIDF" $SERV -D --pid-path=/var/run -R .pid
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

stop(){
 if   log_warning_msg "$NAME stopping"; $SERV -D -F -s -k shutdown
 then log_success_msg "$NAME stopped"
 else log_failure_msg "$NAME not stopped"; status; return 0
 fi
}

reload(){ $SERV -k reconfigure; }
restart(){ stop && sleep 2; start; }
try_restart(){ status && restart; return 0; }
