# trunk/grimoire/metalog.d/init.d.on/metalog
#!/bin/bash

### BEGIN INIT INFO
# Provides: $syslog
# Required-Start: hostname $local_fs run
# Should-Start: hwclock ntpd
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: metalog is a replacement for sysklog and syslog
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/metalog.pid
SERV=/usr/sbin/metalog
OPTS='--sync  --daemonize'

parse "$@"
