# trunk/grimoire/oidentd/xinetd.d/oidentd
service auth
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= daemon
	group			= daemon
	server			= /usr/sbin/oidentd
	server_args		= --stdio
	per_source		= 2
	log_type		= FILE /var/log/xinetd/oidentd
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
