# trunk/grimoire/immune-auto-login.d/init.d/auto-login
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: gui_login
# Required-Start: $network dbus run ICE-unix
# Should-Start: $remote_fs udev xfs dbus kdm
# Required-Stop: kill
# Should-Stop:  $remote_fs initramfs
# Default-Start: 4 5
# Default-Stop: 0 1 2 3 6
# Short-Description: auto-login provides a logged in KDE session
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
server /usr/sbin/auto-login
name auto-login
domain limited
