# trunk/grimoire/dmraid.d/init.d.on/dmraid
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: dmsetup lock sys
# Should-Start: udev
# Required-Stop: halt reboot
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: activates BIOS assisted RAID
### END INIT INFO

. /lib/lsb/init-functions

SERV=/sbin/dmraid
  ON='-ay'
 OFF='-an'

help(){ echo "Usage: $0 {start|stop}" 1>&2; }

case $1 in
 start) OPTS=$ON start ;;
 stop)  OPTS=$OFF stop ;;
 try-restart) : ;;
 *) help ;;
esac
