# trunk/grimoire/apcupsd.d/init.d/UPSdown
#!/bin/sh

### BEGIN INIT INFO
# Default-Mode: 710
# Default-Stop: 0 
# Short-Description: UPSdown shuts down UPS after a computer shutdown
### END INIT INFO

stop(){
 if [ -f /etc/apcupsd/powerfail ]; then
  echo "APCUPSD will now power off the UPS"
  /etc/apcupsd/apccontrol killpower
  echo "Verify the UPS shuts down or turn off the system"
 fi
}

case $1 in
 stop) stop ;;
esac
