# trunk/grimoire/apcupsd.d/init.d/UPSdown
#!/bin/sh

### BEGIN INIT INFO
# Default-Mode: 500
# Default-Stop: 0 
# Short-Description: UPSdown shuts down UPS after a computer shutdown
### END INIT INFO

if [ "$1" == try-restart ]; then exit 2; fi

stop(){
 if [ -f /etc/apcupsd/powerfail ]; then
  echo "APCUPSD will now power off the UPS"
  /etc/apcupsd/apccontrol killpower
  echo "Verify the UPS shuts down or turn off the system"
 fi
}
