# trunk/grimoire/apcupsd.d/init.d/UPSdown
#!/bin/sh

### BEGIN INIT INFO
# Provides: UPSdown
# Required-Start: 
# Required-Stop: 
# Default-Start: 
# Default-Stop: 0 
# Short-Description: UPSdown shuts down UPS after a computer shutdown
# Description: UPSdown shuts down UPS after a computer shutdown during
#     a power failure
### END INIT INFO

if [ -f /etc/apcupsd/powerfail ]; then
        echo
        echo "APCUPSD will now power off the UPS"
        echo
        /etc/apcupsd/apccontrol killpower
        echo
        echo "Verify the UPS shuts down or turn off the system"
        echo
fi
#
# End /etc/init.d/UPSdown

