# trunk/grimoire/logind.d/init.d/logind
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $dev
# Should-Start: cgroup-init-rootfs cgroup-init-root hostname rootname
# Required-Stop:  $local_fs kill
# Default-Start: S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: provides dynamic tty allocation login daemon
### END INIT INFO

[ -f /sbin/logind ] || exit 0

. /lib/lsb/init-functions

only start stop status try-restart
server /sbin/logind
domain user
