# trunk/grimoire/immune-sentinel-memory.d/init.d/sentinel-memory
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: Prevents memory depletion
### END INIT INFO

[ -f /usr/sbin/sentinel-memory ] || exit 0

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

# Modifications to this file are LOST!
# Set variables SDIR and SIZE to alternate values only in
# /etc/init.d/conf.d/sentinel-memory

SDIR=/aux/run/swap		# directory which will hold swap files
SIZE=1024			# size of swap file

. /lib/lsb/init-functions

fast
only start stop restart status configure display
server	/usr/sbin/sentinel-memory
options	"$SDIR ${SIZE}M"
cgroup realtime

display(){
 echo "Swap Directory $SDIR"
 echo "Swap file size $SIZE"
}

configure(){
 input(){ dialog --stdout --no-cancel --timeout 60  --title $"$NAME configuration menu" --inputbox "$1" 10 50 $2; }

 SIZE=$( input $"Megabytes per swap file?"  "${SIZE:=1024}" )
 SDIR=$( input $"Directory for swap files?" "${SDIR:=/aux/run/swap}" )

 if   (( SIZE > 0 )) && [ -n "$SDIR" ]
 then echo "SDIR=$SDIR; SIZE=$SIZE" > /etc/init.d/conf.d/$NAME
 fi
 try_restart
}
