# trunk/grimoire/net-tools.d/init.d.on/loopback
#!/bin/bash

### BEGIN INIT INFO
# Provides: loopback
# Required-Start: $local_fs
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: loopback starts the loopback device
# Description: loopback starts the loopback device
### END INIT INFO

. /lib/lsb/init-functions

start() {
 log_warning_msg "loopback configuring"
 /bin/hostname -F /etc/hostname
 /sbin/ifconfig lo 127.0.0.1
 log_success_msg "loopback configured"
}

case $1 in
 start) start ;;
esac
