# trunk/grimoire/immune-sentinel-sync.d/init.d.on/sentinel-sync
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: synchronizes installed software on boxes weekly
# To modify DAYS set DAYS in /etc/init.d/conf.d/sentinel-sync
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

DAYS=7

. /lib/lsb/init-functions

server  /usr/sbin/sentinel-sync
options $DAYS

parse "$@"
