# trunk/grimoire/immune-sentinel-sync.d/init.d.on/sentinel-sync
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: synchronizes installed software on boxes weekly
# To modify DAYS set DAYS in /etc/init.d/conf.d/sentinel-sync
### END INIT INFO

DAYS=7

. /lib/lsb/init-functions

pid_file /var/run/sentinel-sync.pid
server   /sbin/sentinel-sync
options  $DAYS

parse "$@"
