# trunk/grimoire/immune-sentinel-sync.d/init.d.on/sentinel-sync
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: synchronizes installed software on boxes weekly
# To modify DAYS set DAYS in /etc/init.d/conf.d/sentinel-sync
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

DAYS=7

. /lib/lsb/init-functions

fast
only start stop restart try-restart status configure display
server  /usr/sbin/sentinel-sync
options $DAYS

display(){ echo "Days $DAYS"; }

configure(){
 input(){ dialog --stdout --no-cancel --timeout 60  --title $"$NAME configuration menu" --inputbox "$1" 10 50 $2; }

 DAYS=$( input $"Days between synchronization?" ${DAYS:=1} )

 if   (( DAYS > 0 ))
 then echo "DAYS=$DAYS" > /etc/init.d/conf.d/$NAME
 fi
 try_restart
}
