# trunk/grimoire/mdadm.d/init.d.on/mdadm--assemble--scan
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: mdadm
# Required-Start: dmsetup lock sys
# Should-Start: udev mdadm--auto-detect
# Required-Stop: finale halt reboot reboot-kexec
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mdadm activates the software RAID
# Raid partitions should be type 0xFD
# This script assembles the raid arrays using configuration from
# /etc/mdadm.conf
# Shutting down RAIDS during shutdown is probably not required.
### END INIT INFO

if [ "$1" == try-restart ]; then exit 2; fi

. /lib/lsb/init-functions

SERV=/sbin/mdadm
CONF=/etc/mdadm.conf
  ON='--assemble --scan'
 OFF='--stop --scan'

help(){ echo "Usage: $0 {start|stop}" 1>&2; }

case $1 in
 start) OPTS=$ON ;;
 stop)  OPTS=$OFF ;;
esac
