# trunk/grimoire/mdadm.d/init.d.on/mdadm
#!/bin/bash

### BEGIN INIT INFO
# Provides: mdadm
# Required-Start: dmsetup lock sys
# Should-Start: udev
# Required-Stop: halt reboot
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mdadm activates the software RAID
# Raid partitions should be type 0xFD
# Then, the linux kernel can automatically detect and activate RAIDS
# when told to do during sysinit execution on the initramfs by
# mdadm --auto-detect
### END INIT INFO

. /lib/lsb/init-functions

NAME=mdadm
SERV=/sbin/mdadm

start(){
 if   log_warning_msg "$NAME $OPTS"; $SERV --auto-detect
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

stop(){
 if   log_warning_msg "$NAME $OPTS"; $SERV --stop --scan
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
 start) start ;;
 stop) stop ;;
 try-restart) : ;;
 *) help ;;
esac
