# trunk/grimoire/mdadm.d/init.d.on/mdadm--auto-detect
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: mdadm
# Required-Start: dmsetup lock sys
# Should-Start: udev
# Required-Stop: finale halt reboot reboot-kexec
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: mdadm activates the software RAID
# Raid partitions should be type 0xFD
# Then, the linux kernel can automatically detect and activate RAIDS
# when told to do during sysinit execution on the initramfs by
# mdadm --auto-detect
# Shutting down RAIDS during shutdown is probably not required.
### END INIT INFO

. /lib/lsb/init-functions

only start stop
cgroup none
server /sbin/mdadm

case $1 in
 start) options '--auto-detect' ;;
 stop)  options '--stop --scan' ;;
esac
