# trunk/grimoire/xdm.d/init.d/xdm
#!/bin/bash

### BEGIN INIT INFO
# Provides: gui_login
# Required-Start: $local_fs run ICE-unix
# Should-Start: $network $remote_fs udev xfs
# Required-Stop: kill
# Should-Stop: $remote_fs initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Default-Mode: 500
# Short-Description: xdm provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

fast
affine
only start stop restart status
pid_file /var/run/xdm.pid
server /usr/bin/X11/xdm
config /etc/X11/xdm/xdm-config
domain user
