# trunk/grimoire/slmodem.d/init.d.on/slmodemd
#!/bin/bash

### BEGIN INIT INFO
# Provides: slmodemd
# Required-Start: $local_fs
# Should-Start: modprobe run udev
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: slmodemd starts the winmodem
# Description: slmodemd starts the winmodem
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/slmodemd
PIDF=/var/run/slmodemd.pid
NAME='slmodemd'
CONF='-c USA /dev/slamr0'
MODULE='slamr'

status() {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}


start()  {
  run()  {
    modprobe $MODULE
    start_daemon -b -n -20 $SERV $CONF
  }

  if   log_warning_msg "$NAME starting"; run
  then log_success_msg "$NAME started"
  else log_failure_msg "$NAME not started"
  fi
}

stop()   {
  if   log_warning_msg "$NAME stopping"; killproc -p $PIDF -s 15 $SERV
  then log_success_msg "$NAME stopped"
  else log_failure_msg "$NAME not stopped"
  fi
}

help()   {  echo "Usage: $0 {start|stop|status}" 1>&2;    }

case  $1  in
        start)  start  ;;
         stop)  stop   ;;
       status)  status ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
