# trunk/grimoire/slmodem.d/init.d.on/slmodemd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs
# Should-Start: modprobe run udev
# Required-Stop: $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: slmodemd starts the winmodem
### END INIT INFO

. /lib/lsb/init-functions

server /usr/sbin/slmodemd
pid_file /var/run/slmodemd.pid
name slmodemd
MODULE=slamr
options '-c USA /dev/slamr0'

start(){
 run(){
  modprobe $MODULE
  start_daemon -b $SERV $OPTS
 }

 if   log_warning_msg "$NAME starting"; run
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"
 fi
}
