# trunk/grimoire/procps.d/init.d/sysctl
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc sys
# Should-Start: udev modprobe
# Default-Start: S
# Short-Description: sysctl sets kernel parameters
### END INIT INFO

. /lib/lsb/init-functions

only start
server /sbin/sysctl
CONF=/etc/sysctl.conf
options "-p $CONF"
