# trunk/grimoire/smack-util.d/init.d.on/smack
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-execute: proc
# Default-execute: S
# Short-Description: smack initializes access rules and CIPSO bindings
### END INIT INFO

. /lib/lsb/init-functions

only start
name smackcipso
server /sbin/smackcipso
config /etc/smack/cipso

start(){
 run(){
  [ -f /etc/smack/accesses ] && /sbin/smackload        < /etc/smack/accesses
  [ -f /etc/smack/cipso    ] && /sbin/smackcipso       < /etc/smack/cipso
  [ -f /etc/smack/ambient  ] && cat /etc/smack/ambient > /smack/ambient
  [ -f /etc/smack/nltype   ] && cat /etc/smack/nltype  > /smack/nltype
 }
 if ! [ -f $CONF ]; then
 log_warning_msg "$NAME not starting; $CONF does not exist"; return; fi
 log_warning_msg "$NAME starting"; run
 log_success_msg "$NAME started"
}
