# trunk/grimoire/lm_sensors.d/init.d/fancontrol
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: sensors run
# Required-Stop: syslog kill 
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: sensord is lm_sensors hardware health monitoring daemon
### END INIT INFO

. /lib/lsb/init-functions

fast
only start stop restart try-restart status
#pid_file /var/run/fancontrol.pid
server /usr/sbin/fancontrol
config /etc/fancontrol

start(){
 if ! status &>/dev/null; then
  log_warning_msg "$NAME starting"; setsid $SERV &>/var/log/fancontrol &
  log_success_msg "$NAME started"
 fi
}
