# trunk/grimoire/lm_sensors.d/init.d/sensors_modules
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: modprobe udev
# Required-Stop: syslog kill 
# Default-Start: S
# Short-Description: sensors_modules loads modules specified in /etc/sysconfig/lm_sensors
### END INIT INFO

. /lib/lsb/init-functions

only start
config /etc/sysconfig/lm_sensors

start(){
 load_mods(){
  log_warning_msg "sensor modules loading";
  . $CONF
  for module in $BUS_MODULES $HWMON_MODULES
  do  /sbin/modprobe $module
  done
  log_success_msg "sensor modules loaded"
 }

 if   [ -f $CONF ]
 then load_mods
 else log_warning_msg "Please run sensors-detect to select kernel modules to load."
 fi
}
