# trunk/grimoire/nfs-utils.d/init.d.on/exportfs
#!/bin/bash

### BEGIN INIT INFO
# Provides: exportfs
# Required-Start: rpcbind var_lib_nfs
# Required-Stop:  kill
# Should-Start: $network
# Should-Stop:  rpcbind $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: exportfs 
rpc.nfsd implements the user level part of the NFS service.
# Description: rpc.nfsd merely starts the specified number of kernel threads.
### END INIT INFO

. /lib/lsb/init-functions

SERV=/sbin/exportfs
NAME=exportfs
CONF=/etc/exports

case $@ in
 start) OPTS=-r ;;
  stop) OPTS=-ua ;;
esac

parse "$@"
