# trunk/grimoire/nfs-utils.d/init.d.on/lockd
#!/bin/bash

### BEGIN INIT INFO
# Provides: lockd
# Required-Start: 
# Should-Start: udevd udevsettle
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: lockd loads the lockd module
### END INIT INFO

. /lib/lsb/init-functions

NAME=lockd

parse "$@"
