# trunk/grimoire/nfs-utils.d/init.d.on/lockd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Should-Start: udevd udevsettle
# Default-Start: S
# Short-Description: lockd loads the lockd module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
