# trunk/grimoire/nfs-utils.d/init.d.on/nfs
#!/bin/bash

### BEGIN INIT INFO
# Provides: nfs
# Required-Start: 
# Should-Start: udevd udevsettle
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: nfs loads the nfs module
### END INIT INFO

. /lib/lsb/init-functions

NAME=nfs

parse "$@"
