# trunk/grimoire/nfs-utils.d/init.d.on/nfs
#!/bin/bash

### BEGIN INIT INFO
# Should-Start: udevd udevsettle
# Default-Start: S
# Short-Description: nfs loads the nfs module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
