# trunk/grimoire/nfs-utils.d/init.d.on/nfsd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Should-Start: udevd udevdsettle modules
# Default-Start: S 2 3 4 5
# Short-Description: nfsd mounts nfsd on /proc/fs/nfsd
# Description: nfsd mounts nfsd on /proc/fs/nfsd
#	This filesystem needs to be mounted before most
#	NFS daemons, particularly exportfs, mountd, svcgssd, idmapd.
### END INIT INFO

. /lib/lsb/init-functions
only start stop restart status
name /proc/fs/nfsd
config /etc/exports
options '-n -t nfsd nfsd'
