# trunk/grimoire/nfs-utils.d/init.d.on/var_lib_nfs
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: proc
# Default-Start: 2 3 4 5
# Short-Description: var_lib_nfs mounts tmpfs on /var/lib/nfs
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /var/lib/nfs
options '-o defaults,size=1m,nr_inodes=256,mode=0755 -t tmpfs tmpfs'
