# trunk/grimoire/sysvinit.d/init.d/halt
#!/bin/bash

### BEGIN INIT INFO
# Provides: halt
# Required-Start: 
# Required-Stop: 
# Default-Start: 
# Default-Stop: 0
# Short-Description: halt kills lingering processes and powers down the box
# Description: halt kills lingering processes and powers down the box
### END INIT INFO

. /lib/lsb/init-functions

log_success_msg "Halting"
halt  -d -f -i -p
