# trunk/grimoire/sysvinit.d/init.d/reboot
#!/bin/bash

### BEGIN INIT INFO
# Provides: reboot
# Required-Start: 
# Required-Stop: 
# Default-Start: 
# Default-Stop: 6
# Short-Description: reboot kills lingering processes and reboots the box
# Description: reboot kills lingering processes and reboots the box
### END INIT INFO

. /lib/lsb/init-functions

log_success_msg "Rebooting"
reboot  -d -f -i
