# trunk/grimoire/heimdal.d/init.d/hpropd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: $network run
# Required-Stop: $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: hpropd recieves databases sent by hprop, use ipropd-slave instead
# Description: hpropd recieves databases sent by hprop
# hpropd recieves databases sent by hprop
# and writes it as a local database.
# However, the combination of using ipropd-slave and ipropd-master
# to send incremental updates may be better,
# especially when the database is large.
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart reload status
pid_file /var/run/hpropd.pid
server /usr/libexec/hpropd
name hpropd
options --no-inetd

reload(){
 PIDS="$( pidofproc -p $PIDF $SERV )"
 if   [ -n            "$PIDS" ]
 then kill  -s SIGHUP  $PIDS  &&
      log_success_msg "reload"
 fi
}
