# trunk/grimoire/heimdal.d/init.d/hpropd
#!/bin/bash

### BEGIN INIT INFO
# Provides: hpropd
# Required-Start: $network run
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: hpropd recieves databases sent by hprop, use ipropd-slave instead
# Description: hpropd recieves databases sent by hprop
# hpropd recieves databases sent by hprop
# and writes it as a local database.
# However, the combination of using ipropd-slave and ipropd-master
# to send incremental updates may be better,
# especially when the database is large.
### END INIT INFO

. /lib/lsb/init-functions

PIDF="/var/run/hpropd.pid"
SERV="/usr/libexec/hpropd"
NAME="hpropd"

status() {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}

start() {
  log_warning_msg  "$NAME starting"
  if    start_daemon     -p $PIDF $SERV --no-inetd
  then  log_success_msg  "$NAME started"
  else  log_failure_msg  "$NAME not started"
        status
  fi
}

stop() {
  log_warning_msg "$NAME stopping"
  if    killproc        -p $PIDF -s 15  $SERV
  then  log_success_msg "$NAME stopped"
  else  log_failure_msg "$NAME not stopped"
        status
        return 0
  fi
}

reload() {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if    [ -n            "$PIDS" ]
  then  kill  -s SIGHUP  $PIDS  &&
        log_success_msg "$NAME configuration reloaded"
  fi
}

restart() {  stop  &&  sleep  1;  start;  }

try_restart() {
  status  &&  restart
  return 0
}

help() {
  echo	"Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}

case $1 in
       restart)      restart  ;;
   try-restart)  try_restart  ;;
        reload)  reload       ;;
  force-reload)  reload       ;;
         start)  start        ;;
        status)  status       ;;
          stop)  stop         ;;
             *)  help         ;;
esac
