# trunk/grimoire/heimdal.d/init.d/kpasswdd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: $network run
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: kpasswdd serves requests for password changes.
# Description: kpasswdd serves requests for password changes.
# It listens on UDP port 464 (service kpasswd) and
# processes requests when they arrive.
# It changes the database directly and
# should thus only run on the master KDC.
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/kpasswdd.pid
SERV=/usr/libexec/kpasswdd
NAME=kpasswdd

reload(){
 PIDS="$( pidofproc -p $PIDF $SERV )"
 if   [ -n            "$PIDS" ]
 then kill  -s SIGHUP  $PIDS  &&
      log_success_msg "$NAME configuration reloaded"
 fi
}

parse "$@"
