# trunk/grimoire/immune-boot-blaze.d/init.d.on/boot-blaze
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs
# Should-Start: tmp
# Default-Start: 2 3 4 5
# Short-Description: boot-blaze reduces system startup and login time
### END INIT INFO

. /lib/lsb/init-functions

server /usr/sbin/bootblaze
parse "$@"
