# trunk/grimoire/immune-boot-blaze.d/init.d.on/boot-blaze
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Should-Start: tmp
# Default-Start: 2 3 4 5
# Short-Description: boot-blaze reduces system startup and login time
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

. /lib/lsb/init-functions

fast
cgroup idle
only start reset
server /usr/sbin/boot-blaze

reset(){ rm -f /etc/boot-blaze; }
