# trunk/grimoire/cups.d/init.d/cups
#!/bin/bash

### BEGIN INIT INFO
# Provides: cups
# Required-Start: $network run
# Required-Stop: $kill
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: cups is a print spooler
# Description: cups provides a portable printing layer
### END INIT INFO

. /lib/lsb/init-functions

#PIDF=/var/run/cupsd.pid
SERV=/usr/sbin/cupsd
NAME=cupsd

parse "$@"
