# trunk/grimoire/heimdal.d/xinetd.d/ftp
# This file is customizable.
# For instructions run:
# man xinetd.conf

service ftp
{
	disable		= no
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/libexec/ftpd
	per_source	= 8
	log_type	= FILE /var/log/xinetd/ftp
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
