# trunk/grimoire/heimdal.d/xinetd.d/kerberos-adm
service kerberos-adm
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/libexec/kadmind
	per_source	= 8
	log_type	= FILE /var/log/xinetd/kerberos-adm
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
