# trunk/grimoire/gdm.d/init.d/gdm
#!/bin/bash

### BEGIN INIT INFO
# Provides: gdm gui_login
# Required-Start: $network dbus run ICE-unix
# Should-Start: $remote_fs udev xfs dbus
# Required-Stop: kill
# Should-Stop: $remote_fs initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: gdm provides an x11 login.
# Description: gdm provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/gdm.pid
SERV=/usr/sbin/gdm
NAME=gdm

parse "$@"
