# trunk/grimoire/heimdal/xinetd.d/ekshell
service ekshell
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/libexec/rshd
        server_args		= -kx
	per_source		= 4
	log_type		= FILE /var/log/xinetd/ekshell
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
