# trunk/grimoire/martian.d/init.d.on/martian
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs martian_dev
# Should-Start: udev
# Required-Stop: $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: martian starts the winmodem
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/martian_modem
NAME='martian linmodem'
OPTS='--daemon --syslog --realtime'

parse "$@"
