# trunk/grimoire/martian.d/init.d.on/martian
#!/bin/bash

### BEGIN INIT INFO
# Provides: martian
# Required-Start: $local_fs
# Should-Start: udev modprobe
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: martian starts the winmodem
# Description: martian starts the winmodem
### END INIT INFO

. /lib/lsb/init-functions

SERV=/usr/sbin/martian_modem
NAME='martian linmodem'

start() {
 run() {
  modprobe martian_dev
  start_daemon $SERV --daemon --syslog --realtime
 }

 if   log_warning_msg "$NAME starting"; run
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"
 fi
}

parse
