# trunk/grimoire/martian.d/init.d.on/martian_dev
#!/bin/bash

### BEGIN INIT INFO
# Provides: martian_dev
# Required-Start: $local_fs
# Should-Start: udev
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Default-Stop:
# Short-Description: martian_dev loads the module
# Description: martian_dev loads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME=martian_dev

parse "$@"
