# trunk/grimoire/module-init-tools.d/init.d/depmod
#!/bin/bash

### BEGIN INIT INFO
# Provides: depmod
# Required-Start: $local_fs
# Should-Start:
# Required-Stop:  
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: depmod creates a list of module dependencies
# Description: depmod creates a list of module dependencies
### END INIT INFO

. /lib/lsb/init-functions

start()  {
if  [  -x  /sbin/depmod  ];  then
    log_warning_msg  "depmod -Aa about to run"
    depmod  -Aa
    log_success_msg  "depmod -Aa complete"
  fi
}

help()  {
  echo "Usage: $0 {start}" 1>&2
}

case  $1  in
  start)  start  ;;
      *)  help   ;;
esac
