# trunk/grimoire/nfs-utils.d/init.d/lockd
#!/bin/bash

### BEGIN INIT INFO
# Provides: lockd
# Required-Start: rpc.statd
# Required-Stop:  rpc.statd kill
# Should-Start: portmap $network
# Should-Stop:  portmap $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: lockd starts the linux kernel module lockd
# Description: lockd starts the linux kernel module lockd
### END INIT INFO

. /lib/lsb/init-functions

NAME=lockd
#OPTS=

start(){
 if   log_warning_msg  "$NAME starting"; modprobe lockd $OPTS
 then log_success_msg  "$NAME started"
 else log_failure_msg  "$NAME not started"; status
 fi
}

stop(){
 log_warning_msg "$NAME stopping"; rmmod lockd
 log_success_msg "$NAME stopped"
}

parse "$@"
