# trunk/grimoire/nfs-utils.d/init.d/nfsd
#!/bin/bash

### BEGIN INIT INFO
# Provides: nfsd
# Required-Start: proc
# Should-Start: 
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: nfsd mounts nfsd on /proc/fs/nfsd
# Description: nfsd mounts nfsd on /proc/fs/nfsd
#	This filesystem needs to be mounted before most
#	NFS daemons, particularly exportfs, mountd, svcgssd, idmapd.
### END INIT INFO

. /lib/lsb/init-functions

NAME='/proc/fs/nfsd'
OPTS='-n -t nfsd nfsd'

parse "$@"

