# trunk/grimoire/samba.d/init.d.on/smbd
#!/bin/bash

### BEGIN INIT INFO
# Provides: smbd
# Required-Start: nmbd run
# Should-Start:
# Required-Stop: nmbd
# Should-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: smbd provides a SMB/CIFS service to clients
# Description: smbd provides a SMB/CIFS service to clients
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/smbd.pid
SERV=/usr/sbin/smbd
CONF=/etc/samba/smb.conf
NAME=smbd
OPTS=-D

reload(){
 PIDS="$( pidofproc -p $PIDF $SERV )"
 if   [ -n            "$PIDS" ]
 then kill  -s SIGHUP  $PIDS  &&
      log_success_msg "$NAME configuration reloaded"
 fi
}

help(){
 echo "Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}

case $1 in
       reload) reload ;;
 force-reload) reload ;;
*) parse "$@" ;;
esac
