# trunk/grimoire/kde-workspace.d/init.d/kdm
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: gui_login
# Required-Start: $local_fs run var_tmp ICE-unix var_lib_kdm
# Should-Start: $network $remote_fs udev xfs hald
# Required-Stop: kill
# Should-Stop: $remote_fs $network initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: kdm provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

fast
affine
only start stop restart status
server /usr/bin/kdm
config /usr/share/config/kdm/kdmrc
domain user
cgroup fast

restart(){ stop; sleep 45; start; }
