# trunk/grimoire/kde-workspace.d/init.d/var_lib_kdm
#!/bin/bash

# Copyright Kyle Sallee 2013
# All rights reserved.
# For use with the distribution sorcerer only!

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 4 5
# Short-Description: creates /var/lib/kdm/ and assigns proper mode and ownership
### END INIT INFO

. /lib/lsb/init-functions

only start
name /var/lib/kdm

chkdir(){
 if    ! [ -d /var/lib/kdm ]; then
  if   log_warning_msg "rm -f /var/lib/kdm && mkdir -pm 755 /var/lib/kdm/" &&
                        rm -f /var/lib/kdm && mkdir -pm 755 /var/lib/kdm/
  then log_success_msg "/var/lib/kdm/ created"
  else log_failure_msg "/var/lib/kdm/ creation failed"; return 1
  fi
 elif    [ -d /var/lib/kdm ] &&
       ! [ -n "$( find /var/lib/kdm -maxdepth 0 -uid 0 -gid 0 -perm 755 )" ]
 then
  if   log_warning_msg "chown root:root /var/lib/kdm/ && chmod 755 /var/lib/kdm/" &&
                        chown root:root /var/lib/kdm/ && chmod 755 /var/lib/kdm/
  then log_success_msg "/var/lib/kdm/ ownership and mode adjusted"
  else log_failure_msg "/var/lib/kdm/ ownership and mode adjustment failed"; return 1
  fi
 fi
}

start(){
 if   log_warning_msg "/var/lib/kdm/ existence ownership mode checking"; chkdir
 then log_success_msg "/var/lib/kdm/ verified"
 else log_failure_msg "/var/lib/kdm/ improper"
 fi
}
