# trunk/grimoire/smartmontools.d/init.d.on/smartd
#!/bin/bash

### BEGIN INIT INFO
# Provides: smartd
# Required-Start: $syslog modprobe run
# Should-Start: udev
# Required-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop:  0 6
# Short-Description: smartd monitors health status of hard drives
# Description: smartd monitors health status of hard drives
### END INIT INFO

. /lib/lsb/init-functions

NAME=smartd
PIDF=/var/run/smartd.pid
SERV=/usr/sbin/smartd
OPTS="-p $PIDF"

parse "$@"
