# trunk/grimoire/smartmontools.d/init.d.on/smartd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $syslog modprobe run
# Should-Start: udev
# Required-Stop: kill
# Default-Start: 2 3 4 5
# Default-Stop:  0 6
# Short-Description: smartd monitors health status of hard drives
### END INIT INFO

. /lib/lsb/init-functions

fast
only start stop restart try-restart status
pid_file /var/run/smartd.pid
server /sbin/smartd
options "-p $PIDF"
