# trunk/grimoire/hc-cron.d/init.d.on/hc-cron
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Provides: $cron
# Required-Start: $local_fs run
# Required-Stop: kill
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: hc-cron is a vixie cron replacement
### END INIT INFO

. /lib/lsb/init-functions

export TERM=dumb
PIDF=/var/run/crond.pid
SERV=/usr/sbin/crond
CONF=/var/spool/cron/

parse "$@"
