# trunk/grimoire/hc-cron.d/init.d.on/hc-cron
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Provides: $cron
# Required-Start: $local_fs run
# Required-Stop: kill
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: hc-cron is a vixie cron replacement
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
export TERM=dumb
pid_file /var/run/crond.pid
server /usr/sbin/crond
config /var/spool/cron/
