# trunk/grimoire/dbus.d/init.d.on/dbus
#!/bin/bash

### BEGIN INIT INFO
# Provides: dbus
# Required-Start: $local_fs run
# Should-Start:
# Required-Stop: kill
# Should-Stop: kill
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: dbus is a systemwide message bus
# Description: dbus is a systemwide message bus
### END INIT INFO


. /lib/lsb/init-functions

PIDF=/var/run/dbus/pid
SERV=/usr/bin/dbus-daemon
NAME=dbus-daemon
MAID=/var/lib/dbus/machine-id
 VLS=/var/lock/subsys
VLSM=/var/lock/subsys/messagebus

status()  {
  PIDS="$( pidofproc -p $PIDF $SERV )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}


start()  {
  local  STATUS
  log_warning_msg  "$NAME starting"
  status  &>/dev/null
  STATUS=$?
  if    (( STAUTS == 3 ))  &&  [  -f  $MAID  ];  then  rm  -f  $MAID
  elif  (( STATUS == 1 ));                       then  rm  -f  $PIDF
  elif  (( STATUS == 0 ));                       then  status; exit 1
  fi
  [  -d  /var/run/dbus             ]  ||  mkdir  -p  /var/run/dbus
  [  -d  /var/lib/dbus             ]  ||  mkdir  -p  /var/lib/dbus
  [  -d  $VLS                      ]  ||  mkdir  -p  $VLS
  if    start_daemon  -p $PIDF $SERV --system
  then  dbus-uuidgen  --ensure
        echo  -n  ""  >  $VLSM
        log_success_msg  "$NAME started"
  else  log_failure_msg  "$NAME not started"
        status
  fi
}


stop()  {
  log_warning_msg "$NAME stopping"
  if    killproc  -p $PIDF -s 15  $SERV
  then  rm        -f $PIDF $VLSM $MAID
  log_success_msg "$NAME stopped"
  else  log_failure_msg "$NAME not stopped"
        status
        return 0
  fi
}


try_restart() {
  status  &&  restart
  return 0
}


help()  {
  echo "Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}


case  $1  in
       restart)      restart  ;;
   try-restart)  try_restart  ;;
        reload)      restart  ;;
  force-reload)      restart  ;;
         start)  start        ;;
        status)  status       ;;
          stop)  stop         ;;
             *)  help         ;;
esac
