# trunk/grimoire/dbus.d/init.d.on/dbus-uuidgen
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: var_lib_dbus
# Default-Start: 1 2 3 4 5
# Short-Description: creates /var/lib/dbus/machine-id
### END INIT INFO

. /lib/lsb/init-functions

NAME=dbus-uuidgen
OPTS=--ensure
MAID=/var/lib/dbus/machine-id

start(){
 if   log_warning_msg "$NAME $OPTS running"; $NAME $OPTS
 then log_success_msg "$NAME $OPTS ran"
 else log_failure_msg "$NAME $OPTS not run"
 fi
}

stop(){
 if   log_warning_msg "$MAID removing"; rm -f "$MAID"
 then log_warning_msg "$MIAD removed"
 else log_warning_msg "$MAID not removed"
 fi
}

parse "$@"
