# trunk/grimoire/dbus.d/init.d.on/dbus-daemon
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: dbus
# Required-Start: dbus-uuidgen var_run_dbus
# Required-Stop: kill
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: systemwide message bus
### END INIT INFO

. /lib/lsb/init-functions

NAME=dbus-daemon
SERV=/usr/bin/dbus-daemon
OPTS=--system
PIDF=/var/run/dbus/pid

parse "$@"
