# trunk/grimoire/dbus.d/init.d.on/var_lib_dbus
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs
# Default-Start: 1 2 3 4 5
# Short-Description: creates directory /var/lib/dbus
### END INIT INFO

. /lib/lsb/init-functions

only start
name /var/lib/dbus

start(){
 run(){ mkdir -pm 775 $NAME && chown messagebus:messagebus $NAME; }

 if   [ -d $NAME ]; then return; fi
 if   log_warning_msg "$NAME creating"; run
 then log_success_msg "$NAME created"
 else log_failure_msg "$NAME not created"; return 1
 fi
}
