# trunk/grimoire/coda.d/init.d/coda
#!/bin/bash

### BEGIN INIT INFO
# Provides: coda
# Required-Start:
# Required-Stop:
# Should-Start: udev
# Default-Start: 2 3 4 5
# Default-Stop:
# Short-Description: coda loads the module
# Description: coda loads the module
### END INIT INFO

. /lib/lsb/init-functions

NAME=coda

parse "$@"
