# trunk/grimoire/coda.d/init.d/coda
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Should-Start: udev
# Default-Start: 2 3 4 5
# Short-Description: coda loads the module
### END INIT INFO

. /lib/lsb/init-functions
parse "$@"
