# trunk/grimoire/coda.d/init.d/coda-server-logrotate
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network
# Required-Stop:  $network
# Default-Start: 2 3 4 5
# Short-Description: rotates coda server logs
### END INIT INFO

. /lib/lsb/init-functions

only start
server /usr/sbin/coda-server-logrotate
