# trunk/grimoire/coda.d/init.d/venus.cache
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network
# Required-Stop:  $network
# Should-Start: venus_tmpfs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: venus.cache mounts a 1G tmpfs on /usr/coda/venus.cache
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
name /usr/coda/venus.cache
options '-o defaults,size=1g,nr_inodes=1m,mode=0750 -t tmpfs tmpfs'
