# trunk/grimoire/coda.d/init.d/venus
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network coda vutil
# Required-Stop:  $network
# Should-Start: venus.cache
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: venus provides the coda client.
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart status
server /usr/sbin/venus
name venus
config /etc/coda/venus.conf
