# trunk/grimoire/pulseaudio.d/init.d/pulseaudio
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 400
# Required-Start: cgroup-init $network
# Required-Stop: $local_fs
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: For running pulseaudio in system mode.
#	http://www.freedesktop.org/wiki/Software/PulseAudio/Documentation/User/WhatIsWrongWithSystemWide
#	Running in system wide mode as user pulse is not advised.
### END INIT INFO

. /lib/lsb/init-functions

fast
only start stop restart try-restart status
server /usr/bin/pulseaudio
config /etc/pulse/daemon.conf
options --daemonize --system
