# trunk/grimoire/jackit.d/jackd
#!/bin/bash

### BEGIN INIT INFO
# Provides: jackd
# Required-Start: $local_fs shm run
# Required-Stop:  $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: jackd starts the jack audio server
# Description: jackd starts and stops the jack audio server
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/jackd.pid
SERV=/usr/bin/jackd
NAME=jackd
OPTS='-d alsa -M -s -p 4096'

start(){
 if   log_warning_msg "$NAME starting"; start_daemon -n -5 -b -p "$PIDF" $SERV $OPTS &>/var/log/jackd.log
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"; status
 fi
}

parse "$@"
