# trunk/grimoire/sorcerer-installer.d/init.d/alt_boot
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: import_archives
# Should-Start: cryptsetup mdadm udevsettle vgchange
# Default-Start: S
# Short-Description: boot presents additional boot options when booting from I/R
### END INIT INFO

# Copyright 2007, 2008, 2009, 2010, 2011, 2012 by Kyle Sallee,
# All rights reserved.
# For use with the Sorcerer distribution only.

. /lib/lsb/init-functions

only start

start(){
 block(){ find /dev -type b -printf '%p	%Y\n' | sed '/loop/d'; }

 root_menu(){
  dialog --backtitle alt_boot --stdout --no-cancel --title "Root Device Menu" \
         --menu "Please select device for /" 0 0 0 $( block )
 }

 root(){ echo "$( root_menu ) /media/root auto defaults 0 0" >> /etc/fstab; }

 main_menu(){
  dialog --backtitle alt_boot --item-help --stdout --no-cancel --timeout 60 \
         --title "Install/Rescue Menu" \
         --menu  "Please select an option within 60 seconds" 0 0 0 \
         "Install" "for deployment" \
                   "also can be used for manual install or rescue" \
         "Rescue"  "to boot a previously deployed box" \
                   "mounts a selected file system as root"
 }

 main(){
  case $( main_menu ) in
   Rescue) root ;;
  esac
  redraw
 }

 grep -q boot=IR /proc/cmdline || return

 if   log_warning_msg 'boot options menu'; main
 then log_success_msg 'boot options menu complete'
 else log_failure_msg 'boot options menu incomplete'
 fi
}
