# trunk/grimoire/sorcerer-installer.d/init.d/import_boot
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: import_archives
# Should-Start: udev loadkeys
# Default-Start: S
# Short-Description: archives makes the initramfs more usable.
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 if   [ -d /boot ]
 then log_success_msg '/boot already copied?'
 elif log_warning_msg '/boot copying'; cp -a /media/IR/boot /
 then log_success_msg '/boot copied'
 else log_failure_msg '/boot copying failed'; false
 fi
}
