# trunk/grimoire/sorcerer-installer.d/init.d/ir_find
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: udev loadkeys
# Default-Start: S
# Short-Description: ir_find locates and makes the Sorcerer IR / in /etc/fstab
### END INIT INFO

. /lib/lsb/init-functions

only start

start(){
 block(){ while read; do if [ -b "$REPLY" ]; then echo "$REPLY"; fi; done; }
   try(){ while read; do mtest   "$REPLY" && return; done; false; }

 try_mount(){
  local CX=12
  while (( CX > 0 )); do
   (( CX-- ))
   mount -r "$1" /media/IR && return
   umount        /media/IR
   sleep 5
   false
  done
 }

 mtest(){
  if try_mount "$1"; then
   [ -d   /media/IR/boot/isolinux ] && return
   umount /media/IR
  fi
  return 1
 }

 search(){
  mkdir -p /media/IR
  grep -l 1 /sys/block/*/removable |
  sed    's:/sys/block/:/dev/:
                      s:/removable::' |
  grep -v "/dev/fd[0-9]" | block | try
 }

 manual(){
  echo 'Install/Rescue Image not found.'
  echo 'Discover the block device for'
  echo 'the Install/Rescue Image and'
  echo 'mount that device on /media/IR'
  echo 'Sorry for the inconvenience.'
  echo 'Press enter when prompted for password.'
  echo 'When complete type exit to resume'
  sulogin -t 300 /dev/console &>/dev/console 0</dev/console
 }

 grep -q boot=IR /proc/cmdline || return

 if   log_warning_msg 'Install/Rescue finding'; search
 then log_success_msg 'Install/Rescue found'
 else log_failure_msg 'Install/Rescue not found'; manual &>/dev/console
 fi
}
