# trunk/grimoire/mDNSResponder.d/init.d/mdnsd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Should-Start: $network $remote_fs udev xfs
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 6
# Short-Description: mdnsd provides Apple Darwin Multicast DNS / DNS Service Discovery daemon
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/mdnsd.pid
SERV=/usr/sbin/mdnsd

parse "$@"
