# trunk/grimoire/cpufrequtils.d/init.d.on/cpufreq
#!/bin/dash

### BEGIN INIT INFO
# Required-Start: proc sys
# Should-Start: udev
# Default-Start: 1 2 3 4 5
# Short-Description: Loads modules for governors and sets scaling_governor
# Configuration file is /etc/init.d/conf.d/cpufreq
# Define variables in configuration file.
# Examples:
# GOVERNORS="conservative ondemand performance powersave userspace"
# DEFAULT="performance"
# USER="Yes"

# Modules for any or all governors can be loaded.

# However, DEFAULT specifies only the governor
# for all CPU cores to use during boot until the SA or user changes it.
# Providing additional governors are listed in the $GOVERNOR variable
# the invoking cpufreq_$GOVERNOR will switch all cores
# to the specified governor for example:
# cpufreq_performance  # Maximum speed on all cores
# cpufreq_powersave    # Minimum speed on all cores

# If $USER is specified set to anything
# then users also will be able to change the governor
# and CPU core frequencies instead of just root

### END INIT INFO

. /lib/lsb/init-functions

if [ -n "$USER" ]; then
 chmod 666 /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor \
           /sys/devices/system/cpu/cpu*/cpufreq/scaling_setspeed
fi

for GOV in $GOVERNORS; do 
 log_warning_msg "Running modprobe cpufreq_$GOV"
                          modprobe cpufreq_$GOV
done

if                       [ -x /bin/cpufreq_$DEFAULT ]
then log_warning_msg "Running /bin/cpufreq_$DEFAULT"
                              /bin/cpufreq_$DEFAULT
fi
