# trunk/grimoire/util-linux.d/init.d.on/dmesg
#!/bin/bash

### BEGIN INIT INFO
# Provides: dmesg
# Required-Start: $local_fs $syslog
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: dmesg saves the kernel ring buffer to /etc/init.d/log/dmesg
# Description: dmesg saves the kernel ring buffer to /etc/init.d/log/dmesg
### END INIT INFO

. /lib/lsb/init-functions

start() {
  log_warning_msg  "dmesg writing";  /bin/dmesg -c > /etc/init.d/log/dmesg
  log_success_msg  "dmesg written"
}

help()  { echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
       start)  start ;;
 try-restart)  :     ;;
           *)  help  ;;
esac
