# trunk/grimoire/util-linux.d/init.d.on/utmp-wtmp
#!/bin/bash

### BEGIN INIT INFO
# Provides: utmp-wtmp
# Required-Start: $local_fs
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: utmp-wtmp creates login records /var/run/utmp /var/log/wtmp
# Description:       utmp-wtmp creates login records /var/run/utmp /var/log/wtmp
### END INIT INFO

start()  {
  touch            /var/run/utmp  /var/log/wtmp
  chown  root:tty  /var/run/utmp  /var/log/wtmp
  chmod  0664      /var/run/utmp  /var/log/wtmp
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
         stop)  :      ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
