# trunk/grimoire/util-linux.d/init.d.on/utmp
#!/bin/bash

### BEGIN INIT INFO
# Provides: utmp
# Required-Start: $local_fs
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: utmp creates login record /var/run/utmp
# Description:       utmp creates login record /var/run/utmp
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  if    log_warning_msg  "/var/run/utmp creating"
        touch             /var/run/utmp  &&
        chown  root:tty   /var/run/utmp  &&
        chmod  0664       /var/run/utmp
  then  log_success_msg  "/var/run/utmp created"
  else  log_failure_msg  "/var/run/utmp not created"
  fi
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
         stop)  :      ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
