# trunk/grimoire/ypserv.d/init.d/ypserv
#!/bin/bash

### BEGIN INIT INFO
# Provides: ypserv
# Required-Start: portmap
# Required-Stop: portmap
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: ypserv provides a server or slave for serving NIS data
# Description: ypserv is a server that provides NIS data to clients
### END INIT INFO

. /lib/lsb/init-functions

#PIDF=/var/run/ypserv.pid
SERV=/usr/sbin/ypserv
NAME=ypserv

parse "$@"
