# trunk/grimoire/glib.d/init.d.on/glib-compile-schemas
#!/bin/bash
### BEGIN INIT INFO
# Default-Mode: 510
# Required-Stop: $local_fs
# Should-Stop: initramfs tmp
# Default-Stop: 0 6
# Short-Description: updates gschemas.compiled
### END INIT INFO

# Copyright 2011 by Kyle Sallee, all rights reserved.
# for use with Sorcerer only

. /lib/lsb/init-functions

only stop

DIR=/usr/share/glib-2.0/schemas
COM=$DIR/gschemas.compiled

check(){
 if ! [ -f $COM ]; then return; fi
 find $DIR -type f -cnewer $COM | grep -q .
}

update(){ glib-compile-schemas $DIR; }

stop(){
 if   log_warning_msg "$COM checking"; check
 then log_warning_msg "$COM updating"; update
      log_success_msg "$COM updated"
 fi;  log_success_msg "$COM checked"
}
