# trunk/grimoire/slim.d/init.d/slim
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 710
# Provides: gui_login
# Required-Start: $network dbus run ICE-unix
# Should-Start: $remote_fs udev xfs dbus
# Required-Stop: kill
# Should-Stop: $remote_fs initramfs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: slim provides an x11 login.
### END INIT INFO

. /lib/lsb/init-functions

PIDF=/var/run/slim.pid
SERV=/usr/bin/slim

parse "$@"
