# trunk/grimoire/xfs.d/init.d.on/xfs
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs $network run
# Should-Start: $remote_fs udev
# Required-Stop: kill
# Should-Stop: $remote_fs
# Default-Start: 4 5
# Default-Stop: 0 6
# Short-Description: xfs provides an X11 font server.
### END INIT INFO

. /lib/lsb/init-functions

only start stop restart try-restart status
pid_file /var/run/xfs.pid
server /usr/bin/xfs
