# trunk/grimoire/iptables.d/init.d.on/iptables-restore
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Should-Start: loopback $network
# Default-Start: S
# Short-Description: iptables-restore loads previously saved rules
### END INIT INFO

. /lib/lsb/init-functions

name 'iptables firewall rules'
config /etc/iptables.conf

start(){
 [ -s "$CONF" ] &&
 if   log_warning_msg "$NAME loading"; iptables-restore < $CONF
 then log_success_msg "$NAME loading"
 else log_failure_msg "$NAME not loaded"
 fi
}

stop(){ 
 if   log_warning_msg "$NAME flushing"; iptables -F
 then log_success_msg "$NAME flushed"
 else log_failure_msg "$NAME not flushed"
 fi
}
