# trunk/grimoire/udev.d/init.d.on/coldplug
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S
# Short-Description: Allows SA to prefer previous coldplug
### END INIT INFO

. /lib/lsb/init-functions

start(){
 MOD=/etc/coldplug/modules
 DEV=/etc/coldplug/dev
 USE=/tmp/coldplug

 if grep -qi coldplug /proc/cmdline; then
  log_warning_msg "coldplug discovered in /proc/cmdline"
  log_warning_msg "Previously used coldplug /dev and modules request acknowledged"

  if   [ -f $MOD ] && [ -d $DEV ]; then
   echo > $USE
   log_success_msg "File $MOD and directory $DEV will be used"
  else
   log_failure_msg "Unable to comply, because"
   if ! [ -f $MOD ]; then log_failure_msg "$MOD file is missing"; fi
   if ! [ -d $DEV ]; then log_failure_msg "$DEV dir  is missing"; fi
  fi
 elif ! [ -f $USE ]; then
  log_warning_msg "udev with udevadm trigger will handle coldplug events"
  if [ -f $MOD ] && [ -d $DEV ]; then
   log_warning_msg "If this fails or requires longer than 200 seconds"
   log_warning_msg "then please try rebooting with the kernel parameter:"
   log_warning_msg "coldplug"
   log_warning_msg "to reuse previously used coldplug events"
  fi
 elif [ -f $USE ]; then
  log_success_msg "Previously used coldplug /dev and modules will be used"
 fi
}

       stop(){ :; }
    restart(){ :; }
try_restart(){ :; }

parse "$@"
