# trunk/grimoire/udev.d/init.d.on/dev
#!/bin/bash

### BEGIN INIT INFO
# Provides: dev
# Required-Start: proc
# Should-Start:
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: dev mounts tmpfs on /dev
# Description: dev mounts tmpfs on /dev
### END INIT INFO

. /lib/lsb/init-functions

NAME='/dev'
SERV='mount'
OPTS='-o defaults,size=1m,nr_inodes=4k,mode=0755 -t tmpfs tmpfs'

parse "$@"
