# trunk/grimoire/udev.d/init.d.on/dev
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: proc
# Default-Start: S
# Short-Description: dev mounts tmpfs on /dev
### END INIT INFO

. /lib/lsb/init-functions

name /dev

if   grep -q devtmpfs /proc/filesystems
then options '-o defaults,size=1m,nr_inodes=64k,mode=0755 -t devtmpfs devtmpfs'
else options '-o defaults,size=1m,nr_inodes=64k,mode=0755 -t tmpfs tmpfs'
fi
