# trunk/grimoire/udev.d/init.d.on/dev
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: proc
# Default-Start: S
# Short-Description: dev mounts tmpfs on /dev
### END INIT INFO

. /lib/lsb/init-functions

NAME=/dev
if   grep -q devtmpfs /proc/filesystems
then OPTS='-o defaults,size=1m,nr_inodes=64k,mode=0755 -t devtmpfs devtmpfs'
else OPTS='-o defaults,size=1m,nr_inodes=64k,mode=0755 -t tmpfs tmpfs'
fi
parse "$@"
