# trunk/grimoire/udev.d/init.d.on/dev
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: proc
# Default-Start: S
# Short-Description: dev mounts tmpfs on /dev
### END INIT INFO

. /lib/lsb/init-functions

NAME=/dev
OPTS='-o defaults,size=1m,nr_inodes=4k,mode=0755 -t tmpfs tmpfs'

parse "$@"
