# trunk/grimoire/udev.d/init.d.on/dev-cp
#!/bin/bash

### BEGIN INIT INFO
# Provides: dev-cp
# Required-Start: proc dev sys
# Should-Start:
# Required-Stop:  
# Default-Start: S
# S 1 2 3 4 5
# Default-Stop:
# Short-Description: dev-cp copies /lib/udev/dev/* to /dev/
# Description: dev-cp copies /lib/udev/dev/* to /dev/
### END INIT INFO

. /lib/lsb/init-functions

NAME=dev-cp

start(){
 if   log_warning_msg "$NAME copying"; cp -a /lib/udev/dev/* /dev
 then log_success_msg "$NAME copied"
 else log_failure_msg "$NAME not copied"
 fi
 }

stop(){ :; }

parse "$@"
