# trunk/grimoire/udev.d/init.d.on/udev
#!/bin/bash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: proc sys
# Should-Start:
# Required-Stop:  
# Default-Start: S 1 2 3 4 5
# Default-Stop:
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

NAME=udev
SERV=/sbin/udevd
OPTS='--daemon'

start(){
 tmpfs_dev(){
  local OPTIONS='nr_inodes=4096,size=1m,mode=0755'
  if   !  [  -d       /dev ]
  then mkdir -pm 0755 /dev
  fi &&
  mount -o $OPTIONS -t tmpfs tmpfs /dev &&
  if   [   -d  /lib/udev/dev  ]
  then cp  -a  /lib/udev/dev/*  /dev
  fi
 }

 if   status &>/dev/null
 then log_warning_msg "$NAME restarting"; killproc -s 15 "$SERV"
 else log_warning_msg "$NAME starting"; tmpfs_dev
 fi

 start_daemon $SERV $OPTS
 /sbin/udevtrigger
 /sbin/udevsettle --timeout=60

 if   status &>/dev/null
 then log_success_msg "$NAME started"
 else log_failure_msg "$NAME not started"
 fi
}

help() { echo "Usage: $0 {start|status}" 1>&2; }

case $1 in
       start) start  ;;
 try-restart) start  ;;
      status) status ;;
           *) help   ;;
esac
