# trunk/grimoire/udev.d/init.d.on/udevsettle
#!/bin/bash

### BEGIN INIT INFO
# Provides: udevsettle udev
# Required-Start: proc dev sys dev-cp udevd udevtrigger
# Should-Start:
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: udevsettle waits for queued udev events to be handled
# Description: udevsettle waits for queued udev events to be handled
### END INIT INFO

. /lib/lsb/init-functions

NAME=udevsettle
SERV=/sbin/udevsettle
OPTS='--timeout=60'

parse "$@"
