# trunk/grimoire/udev.d/init.d.on/udev
#!/bin/bash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: proc sys
# Should-Start:
# Required-Stop:  
# Default-Start: S 1 2 3 4 5
# Default-Stop:
# Short-Description: udev populates /dev
# Description: udev turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

NAME=udev
PIDF=/var/run/udev.pid
SERV=/sbin/udevd

status() {
# PIDS="$( pidofproc -p $PIDF $SERV )"
  PIDS="$( pidof -s udevd )"
  if   [ -n "$PIDS" ]; then log_warning_msg "$NAME is running"; return 0
  elif [ -f "$PIDF" ]; then log_warning_msg "$NAME is dead and $PIDF exists"; return 1
                       else log_warning_msg "$NAME is not running"; return 3
  fi
}


start()  {
  ramfs_dev()  {
    /bin/mount  -o nr_inodes=4096,size=1m,mode=0755 -t tmpfs tmpfs /dev
    /bin/cp     -a  /lib/udev/dev/*  /dev
  }

  if    status &>/dev/null
  then  log_warning_msg  "$NAME restarting";  kill  -s 15   $PIDS
# then  log_warning_msg  "$NAME restarting";  killproc  -p  $PIDF  -s 15  $SERV
  else  log_warning_msg  "$NAME starting";    ramfs_dev
  fi

  start_daemon  -p  $PIDF  $SERV  --daemon
  ## /var/run not mounted on initramfs when udev runs during sysinit
  # pidof  udevd   >  $PIDF
  /sbin/udevtrigger
  /sbin/udevsettle  --timeout=10

  if    status  &>/dev/null
  then  log_success_msg  "$NAME started"
  else  log_failure_msg  "$NAME not started"
  fi
}

help()  {  echo "Usage: $0 {start|status}" 1>&2;  }

case  $1  in
        start)  start   ;;
  try-restart)  :       ;;
       status)  status  ;;
            *)  help    ;;
esac
