# trunk/grimoire/udev.d/init.d.on/udevd
#!/bin/bash

### BEGIN INIT INFO
# Provides: udevd udev
# Required-Start: dev-cp sys
# Should-Start:
# Required-Stop: $local_fs
# Default-Start: S
# S 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: udevd populates /dev
# Description: udevd turns on dynamic /dev population and module loading
### END INIT INFO

. /lib/lsb/init-functions

NAME=udevd
SERV=/sbin/udevd
OPTS='--daemon'

parse "$@"
