# trunk/grimoire/udev.d/init.d.on/udevsettle
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: udev
# Required-Start: udevd udevtrigger
# Default-Start: S
# Short-Description: udevsettle waits for queued udev events to be handled
### END INIT INFO

. /lib/lsb/init-functions

server /sbin/udevadm
options "settle --timeout 60"
