# trunk/grimoire/udev.d/init.d.on/udevd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Provides: udev
# Required-Start: dev-cp sys
# Required-Stop: $local_fs kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: udevd provides dynamic /dev populates and module loading
### END INIT INFO

. /lib/lsb/init-functions

server /sbin/udevd
options --daemon

stop(){ kill -s 15 $( /bin/pidof udevd ); }
