# trunk/grimoire/udev.d/init.d.on/udevtrigger
#!/bin/bash

### BEGIN INIT INFO
# Provides: udevtrigger
# Required-Start: proc sys dev-cp udevd
# Should-Start:
# Required-Stop:  
# Default-Start: S
# Default-Stop:
# Short-Description: udevtrigger turns on dynamic /dev population and module loading
# Description: udevtrigger requests kernel device uevents.
# Coldplug events are used by udevd to popualte /dev
### END INIT INFO

. /lib/lsb/init-functions

NAME=udevtrigger
SERV=/sbin/udevtrigger

parse "$@"
