# trunk/grimoire/udev.d/init.d.on/udevd-restart
#!/bin/dash

### BEGIN INIT INFO
# Provides: udev
# Required-Start: udevd
# Default-Start: 1 2 3 4 5
# Short-Description: udevd-restart restarts udevd on real root filesystem.
# Description: udevd started in sysinit is unable to populate /dev
#	following the switch to the real root filesystem.
#	Therefore, restarting udevd is necessary.
### END INIT INFO

. /lib/lsb/init-functions

start(){ /etc/init.d/udevd restart; }
stop(){ :; }

parse "$@"
