# trunk/grimoire/udev.d/init.d.on/udevsettle
#!/bin/bash

### BEGIN INIT INFO
# Provides: udevsettle udev
# Required-Start: udevd udevtrigger
# Default-Start: S
# Short-Description: udevsettle waits for queued udev events to be handled
### END INIT INFO

. /lib/lsb/init-functions

NAME=udevsettle
SERV=/sbin/udevsettle
OPTS='--timeout=60'

parse "$@"
