# trunk/grimoire/immune-sentinel-loadavg.d/init.d/sentinel-loadavg
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: 1 2 3 4 5
# Short-Description: causes delayed execution of gcc when loadavg > limit
### END INIT INFO

if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

. /lib/lsb/init-functions

fast
cgroup idle
only start configure display

display(){ echo "Limit configured as $LIMIT"; }

configure(){
 input(){
  dialog --stdout --no-cancel --timeout 60 \
         --title $"$NAME configuration menu" \
         --inputbox "Loadavg limit?" 0 0 $LIMIT
 }

 LIMIT=$( input )

 if   [ -n "$LIMIT" ]
 then echo "LIMIT=$LIMIT" > /etc/init.d/conf.d/$NAME
 fi
}

start(){
 if [ -z "$LIMIT" ]; then
  LIMIT=$( sed "/^processor/p;d" /proc/cpuinfo | wc -l )
  (( LIMIT++ ))
  echo "LIMIT=$LIMIT" > /etc/init.d/conf.d/$NAME
  log_success_msg "$LIMIT set for sentinel-loadavg"
 fi
}
