# trunk/grimoire/bsdftpd-ssl.d/init.d/bsdftpd-ssl
#!/bin/bash

### BEGIN INIT INFO
# Provides: bsdfttpd-ssl
# Required-Start: $network
# Required-Stop:  $network
# Default-Start: 3 4 5
# Default-Stop:  0 1 2 6
# Short-Description: bsdftp-ssl provides a ftp server
# Description:  bsdftp-ssl provides a ftp server
### END INIT INFO

case $1 in
  start|restart)  echo   "$1ing bsdftpd-ssl, ftp daemon."
                  pkill  "^ftpd$"  &&  sleep  5
                  $FTPD
 ;;
 stop)            echo   "$1ping bsdftpd-ssl"
                  pkill  "^ftpd$"  
 ;;



help()  {
  echo "Usage: $0 {restart|try-restart|reload|force-reload|start|status|stop}" 1>&2
}


case  $1  in
       restart)      restart  ;;
   try-restart)  try_restart  ;;
        reload)  reload       ;;
  force-reload)  reload       ;;
         start)  start        ;;
        status)  status       ;;
          stop)  stop         ;;
             *)  help         ;;
esac
