# trunk/grimoire/clamav.d/init.d/clamd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 510
# Required-Start: $local_fs run
# Required-Stop: $local_fs 
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: clamd is the clamav daemon
# Description: clamd is the clamav daemon
### END INIT INFO

. /lib/lsb/init-functions

#PIDF=/var/run/clamd.pid
SERV=/usr/sbin/clamd
OPTS=/etc/clamd.conf

parse "$@"
