# trunk/grimoire/LVM2.d/init.d.on/dmsetup
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: dm-mod
# Should-Start: udev
# Default-Start: S
# Short-Description: dmsetup creates the /dev/mapper/control device
### END INIT INFO

. /lib/lsb/init-functions

only start
options mknodes
server /sbin/dmsetup
