# trunk/grimoire/LVM2.d/init.d.on/lvm
#!/bin/bash

### BEGIN INIT INFO
# Provides: lvm
# Required-Start: 
# Should-Start: hotplug udev
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: lvm makes logical volumes available during boot
# Description: lvm makes logical volumes available during boot
### END INIT INFO

. /lib/lsb/init-functions


start()  {
  log_warning_msg  'Logical volumes activating'
  [  -d  /var  ]  ||
  /sbin/vgscan
  /sbin/vgchange -a y
  log_success_msg  'Logical volumes mounted'
}

help()  {  echo "Usage: $0 {start}" 1>&2;  }

case  $1  in
        start)  start  ;;
  try-restart)  :      ;;
            *)  help   ;;
esac
