# trunk/grimoire/LVM2.d/init.d.on/vgchange
#!/bin/bash

### BEGIN INIT INFO
# Provides: vgchange
# Required-Start: dmsetup lock sys
# Should-Start: mdadm udev vgscan
# Required-Stop: 
# Default-Start: S
# Default-Stop:
# Short-Description: vgchange activates logical volumes
# Description: vgchange activates logical volumes
### END INIT INFO

. /lib/lsb/init-functions

NAME=vgchange
OPTS='--ignorelockingfailure -ay'
SERV=/sbin/vgchange

start(){
 if   log_warning_msg "$NAME $OPTS"; $SERV $OPTS
 then log_success_msg "$NAME success"
 else log_failure_msg "$NAME failure"
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
       start) start ;;
 try-restart) : ;;
           *) help  ;;
esac
