# trunk/grimoire/kbd.d/init.d.on/loadkeys
#!/bin/bash

### BEGIN INIT INFO
# Provides: loadkeys
# Required-Start:
# Should-Start: udev modprobe
# Required-Stop:
# Default-Start: S
# Default-Stop:
# Short-Description: loadkeys loads the keyboard map
# Description:       loadkeys loads the keyboard map
### END INIT INFO

. /lib/lsb/init-functions

start()  {
  keymap_menu()  {
    keymaps()  {
      find  /usr/share/kbd/keymaps -type f         |
      sed   -n  's:.*/::;s:\.map\.gz:	keymap:p'  |
      LC_ALL=C  sort
    }

  TITLE=$"Keymap Selection Menu"
   HELP=$"Please select desired keyboard mapping within 60 seconds."

  /bin/dialog          \
    --stdout           \
    --no-cancel        \
    --timeout  60      \
    --default-item us  \
    --title  "$TITLE"  \
    --menu   "$HELP"   \
    19 60 10           \
    $( keymaps )
}


  log_warning_msg "keymap loading"

  local          KMAP=ask

  if    [  -f             /etc/keymap.conf  ]
  then  read     KMAP  <  /etc/keymap.conf
  fi

  if    [  -z  "$KMAP"           ]  ||
        [      "$KMAP"  ==  ask  ]
  then           KMAP=$( keymap_menu | tee /dev/console )
  fi

  if    [  -n  "$KMAP"  ]  &&  /bin/loadkeys  $KMAP
  then  log_success_msg "keymap loaded"
  else  log_failure_msg "keymap not loaded";  false
  fi
}


help()  {  echo "Usage: $0 {start}" 1>&2;  }


case  $1  in
        start)  start ;;
  try-restart)  :     ;;
            *)  help  ;;
esac
