# trunk/grimoire/kbd.d/init.d.on/loadkeys
#!/bin/bash

### BEGIN INIT INFO
# Should-Start: udev modprobe
# Default-Start: S
# Short-Description: loadkeys loads the keyboard map
### END INIT INFO

. /lib/lsb/init-functions

start(){
 keymap_menu(){
  keymaps(){
   # find /usr/share/keymaps -type f  |
   ls  -R /usr/share/keymaps          |
   sed -n 's:.*/::;s:\.map\.gz:	keymap:p' |
   LC_ALL=C  sort  -u
  }

   HELP=$"Please select desired keyboard mapping within 60 seconds."

  dialog --stdout --no-cancel --timeout 60 --default-item us \
         --title $"Keymap Selection Menu" \
         --menu $"Please select keyboard mapping within 60 seconds." \
         19 60 10 $( keymaps )
 }


 log_warning_msg "keymap loading"

 if [ -z "$KEYMAP" ]; then
           KEYMAP=$( keymap_menu )
  [     -d /etc/init.d/conf.d ] ||
  mkdir -p /etc/init.d/conf.d
  echo "KEYMAP=$KEYMAP" /etc/init.d/conf.d/loadkeys
  chmod 600             /etc/init.d/conf.d/loadkeys
 fi

 if   [ -n "$KEYMAP" ] && loadkeys $KEYMAP
 then log_success_msg "$NAME $KEYMAP keymap loaded"
 else log_failure_msg "$NAME keymap not loaded"; return 1
 fi
}

help(){ echo "Usage: $0 {start}" 1>&2; }

case $1 in
 start) start ;;
 try-restart) : ;;
 *) help ;;
esac
