# trunk/grimoire/ntp.d/init.d/ntpd
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 500
# Required-Start: $network run
# Should-Start: hwclock
# Required-Stop: kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: ntpd synchronizes system time with a server
### END INIT INFO

. /lib/lsb/init-functions

if [ start == "$1" ]; then echo -n >> /var/log/ntp; fi

only start stop status
server /bin/ntpd
config /etc/ntp.conf
options "-g -c $CONF"
domain limited
