# trunk/grimoire/immune-sentinel-elf.d/init.d.on/sentinel-elf
#!/bin/bash

### BEGIN INIT INFO
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: monitors for and repairs /media/root/lib{,64}
### END INIT INFO

if [ "$1" == try-restart ]; then exit 0; fi

PIDF=/var/run/sentinel-elf.pid
SERV=/sbin/sentinel-elf

. /lib/lsb/init-functions

parse "$@"
