# trunk/grimoire/immune-sentinel-elf.d/init.d.on/sentinel-elf
#!/bin/bash

### BEGIN INIT INFO
# Default-Mode: 700
# Required-Start: $local_fs run
# Required-Stop:  $local_fs kill
# Default-Start: S
# Default-Stop: 0 6
# Short-Description: monitors for and repairs /media/root/lib{,64}
### END INIT INFO

if [ "$1" == try-restart ]; then exit 0; fi
if read < /proc/cmdline; [[ $REPLY =~ boot=IR ]]; then exit 0; fi

. /lib/lsb/init-functions

server /sbin/sentinel-elf

parse "$@"
